# sigma/rules/host/privilege_escalation/t1546_008_pre_auth_logon_screen_trigger_winlogon_child_she.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/privilege_escalation.js
title: T1546.008 - Pre-Auth Logon Screen Trigger (winlogon child shell) (winlogon child shell (high-confidence, near-zero FP))
id: 5fa78eb5-d72e-55f5-ac0c-6a01c16fd6a5
status: experimental
description: 'A shell process (cmd.exe, powershell.exe) spawned as a child of winlogon.exe at the logon screen, indicating a sticky-keys or accessibility binary replacement was triggered before authentication, providing unauthenticated SYSTEM-level access Query block: winlogon child shell (high-confidence, near-zero FP).'
references:
    - MITRE ATT&CK T1546.008
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.privilege-escalation
    - attack.t1546.008
    - attack.g0096
    - attack.g0064
    - attack.g0049
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        process.parent.name: winlogon.exe
        process.name:
            - cmd.exe
            - powershell.exe
            - pwsh.exe
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/privilege_escalation.js
    domain: host
    technique: T1546.008
    technique_name: 'Event Triggered Execution: Accessibility Features'
    os: win
    query_block: 1 of 2
    original_query: |-
        process.parent.name: "winlogon.exe"
        AND process.name: ("cmd.exe" OR "powershell.exe" OR "pwsh.exe")
    sysmon: |-
        // Sysmon EID 1 - winlogon spawning a shell (the trigger)
        ParentImage=*\winlogon.exe
        Image=(*\cmd.exe OR *\powershell.exe OR *\pwsh.exe)
        // This should NEVER happen in normal operations.
        // winlogon spawns userinit.exe and LogonUI.exe, not shells.

        // IFEO Debugger set on accessibility binaries
        TargetObject=*\Image File Execution Options\(sethc|utilman|osk|Magnify|Narrator).exe\Debugger*
    powershell: |-
        # Check for accessibility binary replacement / IFEO
        $bins = @('sethc.exe','utilman.exe','osk.exe','Magnify.exe','Narrator.exe')
        foreach ($b in $bins) {
          $path = "$env:SystemRootSystem32$b"
          $hash = (Get-FileHash $path -Algorithm SHA256 -EA SilentlyContinue).Hash
          Write-Host "$b : $hash"

          # Check IFEO Debugger
          $ifeo = Get-ItemProperty "HKLM:SOFTWAREMicrosoftWindows NTCurrentVersionImage File Execution Options$b" -Name Debugger -EA SilentlyContinue
          if ($ifeo) { Write-Host "  IFEO Debugger: $($ifeo.Debugger) (SUSPICIOUS)" }
        }
    registry: |-
        IFEO Debugger keys:
        HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\
          Image File Execution Options\sethc.exe\Debugger
          Image File Execution Options\utilman.exe\Debugger
          Image File Execution Options\osk.exe\Debugger

        Any Debugger value pointing to cmd.exe, powershell.exe,
        or another binary = accessibility backdoor.
    ossdetect: |-
        Sigma:
        - win_proc_creation_winlogon_child_shell.yml
        - win_registry_ifeo_accessibility.yml
        Elastic:
        - Accessibility Features Modification
    tools: |-
        Sysmon (EID 1 parent-child + EID 13 registry)
        Autoruns (IFEO tab)
        EDR behavioral detection
    actors:
        - APT41
        - APT33
        - APT34
    analyst_notes: The winlogon child shell detection catches the TRIGGER of an accessibility backdoor, not just its installation. The installation (binary replacement or IFEO Debugger setting) can be detected via FIM or registry monitoring, but the winlogon-to-shell parent-child event is the definitive proof that the backdoor was actually used. This is particularly valuable for RDP-accessible systems where an attacker can trigger sethc (5x Shift) or utilman (Win+U) at the logon screen without authenticating.
