# sigma/rules/host/privilege_escalation/t1574_011_service_imagepath_failurecommand_redirect_tighte.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/privilege_escalation.js
title: 'T1574.011 - Service ImagePath / FailureCommand Redirect (Tighten: redirected to a user-writable path / interpreter)'
id: c9a2a798-89a1-5e95-a0db-6781a06dea89
status: experimental
description: 'Sysmon EID 13 modifying ImagePath, FailureCommand, or ServiceDll under a Services key by a non-installer process - redirecting a SYSTEM service binary Query block: Tighten: redirected to a user-writable path / interpreter.'
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.privilege-escalation
    - attack.t1574.011
    - attack.g0096
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        winlog.event_id: 13
        registry.path:
            - '*\Services\*\ImagePath'
            - '*\Services\*\ServiceDll'
        registry.data.strings:
            - '*\AppData\*'
            - '*\Temp\*'
            - '*\ProgramData\*'
            - '*\Public\*'
            - '*powershell*'
            - '*cmd.exe*'
            - '*rundll32*'
    condition: selection
falsepositives:
    - Baseline legitimate updaters that rewrite ImagePath during version upgrades to avoid false positives
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/privilege_escalation.js
    domain: host
    technique: T1574.011
    technique_name: 'Hijack Execution Flow: Services Registry Permissions Weakness'
    os: win
    query_block: 2 of 3
    original_query: |-
        winlog.event_id: 13
        AND registry.path: (*\Services\*\ImagePath OR *\Services\*\ServiceDll)
        AND registry.data.strings: (*\AppData\* OR *\Temp\* OR *\ProgramData\* OR *\Public\* OR *powershell* OR *cmd.exe* OR *rundll32*)
    sysmon: |-
        EventID=13 (RegistryValueSet)
        TargetObject matches:
          *\System\CurrentControlSet\Services\*\ImagePath
          OR *\Services\*\Parameters\ServiceDll
          OR *\Services\*\FailureCommand
          OR *\Services\*\FailureActions
        Image NOT in (services.exe, TrustedInstaller.exe,
          msiexec.exe, a known installer/updater)
        Details points at a user-writable path or interpreter:
          *\AppData\* OR *\Temp\* OR *\ProgramData\*
          OR *powershell* OR *cmd.exe* OR *rundll32*

        // Often paired with a subsequent service start (EID 1
        // child of services.exe running the redirected binary).
    powershell: |-
        # Service config redirects by non-installer (Sysmon EID 13)
        Get-WinEvent -FilterHashtable @{
          LogName='Microsoft-Windows-Sysmon/Operational'; ID=13
        } | Where-Object {
          ($_.Properties[4].Value -match 'Services\\[^\\]+\\(ImagePath|ServiceDll|FailureCommand|FailureActions)') -and
          ($_.Properties[3].Value -notmatch 'services\.exe|TrustedInstaller|msiexec')
        } | Select TimeCreated,
          @{n='Target';e={$_.Properties[4].Value}},
          @{n='NewValue';e={$_.Properties[5].Value}},
          @{n='By';e={($_.Properties[3].Value -split '\\')[-1]}}

        # Audit service key ACLs for weak (writable-by-user) perms
        Get-ChildItem 'HKLM:\SYSTEM\CurrentControlSet\Services' |
          ForEach-Object {
            $acl = Get-Acl $_.PSPath -EA SilentlyContinue
            $weak = $acl.Access | Where-Object {
              ($_.IdentityReference -match 'Users|Everyone|Authenticated Users|INTERACTIVE') -and
              ($_.RegistryRights -match 'FullControl|SetValue|WriteKey') -and
              ($_.AccessControlType -eq 'Allow')
            }
            if ($weak) { [pscustomobject]@{ Service=$_.PSChildName; WeakACE=$weak.IdentityReference } }
          }
    registry: |-
        Service configuration values that grant code execution:
        HKLM\SYSTEM\CurrentControlSet\Services\<svc>\
          ImagePath        - the binary the SCM launches (SYSTEM)
          ServiceDll       - (under \Parameters) for svchost svcs
          FailureCommand   - command run on service failure
          FailureActions   - can specify a run-program action
          ObjectName       - the account the service runs as
                             (LocalSystem = SYSTEM)

        Weak-ACL escalation path:
        - If a low-priv user has SetValue/WriteKey on the service
          key, they rewrite ImagePath to their payload, then
          trigger a service restart (or wait for reboot) - the
          SCM launches their binary as the service account
          (commonly LocalSystem = SYSTEM).
        - FailureCommand/FailureActions variant: set a failure
          command, then crash the service to trigger it.

        Investigation pivots:
        - Compare ImagePath against the service's installer
          provenance; a recently-changed ImagePath on a SYSTEM
          service pointing at AppData/Temp is the signature
        - Enumerate service-key ACLs for non-admin write access
          (the precondition) - tools: accesschk -kvuqsw, PowerUp
    ossdetect: |-
        Sigma:
        - registry_set_service_imagepath_change_susp.yml
        - registry_set_service_dll_susp_path.yml
        - registry_set_service_failure_command.yml

        Atomic Red Team:
        - T1574.011 (service registry permissions tests)
        - T1543.003 (overlapping service-config tests)

        Hayabusa:
        - ServiceImagePathModification, ServiceDllHijack rules

        Velociraptor:
        - Windows.System.Services (config + ACL enumeration)
        - Windows.Detection.ServiceRegistryHijack
    tools: |-
        PowerUp (Get-ModifiableService / Invoke-ServiceAbuse)
        SharpUp - finds modifiable services / weak service ACLs
        accesschk.exe -kvuqsw (Sysinternals) - service key ACL audit
        WinPEAS / Seatbelt - enumerate weak service permissions
        Metasploit - service_permissions LPE module
        sc.exe config <svc> binPath= ... - the built-in redirect
        reg.exe - direct ImagePath rewrite where key is writable
    actors:
        - APT41
    analyst_notes: 'This is one of the classic Windows local-privilege-escalation primitives and it''s distinct from the persistence framing of services: here the goal is escalation via a weak service-key ACL. The precondition is that a non-admin user holds SetValue/WriteKey on a service''s registry key (or has write access to the service binary itself - see the file-permissions row). When that''s true, the attacker rewrites ImagePath (or ServiceDll for svchost-hosted services, or FailureCommand) to point at their payload, then triggers a restart - and the Service Control Manager launches it as the service account, almost always LocalSystem. The detection that matters most: Sysmon EID 13 writes to a Services\*\ImagePath|ServiceDll|FailureCommand value by anything other than services.exe/TrustedInstaller/msiexec/a known installer, especially when the new value points at AppData/Temp or an interpreter. Correlate with a subsequent service-start (7036 ''running'' or an EID 1 child of services.exe). The proactive hardening pivot is to enumerate service-key ACLs for non-admin write access (accesschk -kvuqsw / PowerUp) and fix them - that closes the technique entirely. Baseline legitimate updaters that rewrite ImagePath during version upgrades to avoid false positives.'
