# sigma/rules/host/defense_evasion/t1036_004_masquerade_task_or_service_systemd_unit_cron_nam_2.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/defense_evasion.js
title: T1036.004 - Masquerade Task or Service (systemd unit / cron named to blend) (Service spawning a child from a writable pa
id: c0c450eb-9ce9-51af-9639-56ffcd6bcc64
status: experimental
description: 'A persistence service, timer, or cron job given a legitimate-sounding name (systemd-update, dbus-org-helper, network-manager-dispatch, kworker-timer) to evade review of the service and scheduled-task inventory, while its ExecStart or command points to a payload in a writable or non-standard path Query block: Service spawning a child from a writable path (runtime tell): a "systemd-*" or system-sounding unit whose process executes /tmp, /dev/shm, /home, /var/tmp.'
references:
    - MITRE ATT&CK T1036.004
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.defense-evasion
    - attack.t1036.004
    - attack.g0096
    - attack.g0094
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        process.parent.name: systemd
        process.executable:
            - '*/tmp/*'
            - '*/dev/shm/*'
            - '*/home/*'
            - '*/var/tmp/*'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/defense_evasion.js
    domain: host
    technique: T1036.004
    technique_name: 'Masquerading: Masquerade Task or Service'
    os: linux
    query_block: 3 of 3
    original_query: |-
        process.parent.name: ("systemd")
        AND process.executable: (*/tmp/* OR */dev/shm/* OR */home/* OR */var/tmp/*)
    sysmon: |-
        // Sysmon for Linux EID 1 - task/service masquerade
        // (The artifact is a unit/cron file; the masquerade is the NAME.
        //  Detection = legit-sounding name + illegitimate ExecStart/owner.
        //  Cross-ref persistence T1543.002 (systemd) / T1053.003 (cron).)

        // New/modified systemd unit with a system-ish name
        Image=(*/cp OR */mv OR */tee OR */vi OR */vim OR */nano OR */systemctl)
          target in /etc/systemd/system /usr/lib/systemd/system
                    ~/.config/systemd/user /run/systemd/system
          with a name like systemd-* dbus-* network-* udev-* (mimicry)

        // New cron entry with an innocuous name but odd payload
        Image=*/crontab CommandLine matches: *-e*  (interactive edit)
        target in /etc/cron.d /etc/cron.daily /var/spool/cron

        // Auditd
        -w /etc/systemd/system -p wa -k systemd_unit
        -w /usr/lib/systemd/system -p wa -k systemd_unit
        -w /etc/cron.d -p wa -k cron_change
        -w /etc/cron.daily -p wa -k cron_change
        -w /var/spool/cron -p wa -k cron_change
    powershell: |-
        #!/bin/bash
        # T1036.004 - Task/service masquerade hunt
        # A service NAME can be anything; judge by its ExecStart path and
        # package ownership, not by how official the name sounds.

        echo "[*] === systemd units with ExecStart in writable/odd paths ==="
        for u in /etc/systemd/system/*.service /usr/lib/systemd/system/*.service \
                 /run/systemd/system/*.service ~/.config/systemd/user/*.service; do
          [ -f "$u" ] || continue
          es=$(grep -E "^ExecStart=" "$u" 2>/dev/null | head -1 | sed 's/^ExecStart=//;s/^[@+!-]*//')
          binpath=$(echo "$es" | awk '{print $1}')
          case "$binpath" in
            /tmp/*|/dev/shm/*|/home/*|/var/tmp/*|/root/*|*/.*)
              echo "[FLAG] $u"
              echo "    name looks like: $(basename "$u")"
              echo "    ExecStart -> $es" ;;
          esac
        done

        echo ""
        echo "[*] === systemd unit FILES not owned by a package ==="
        HAVE_RPM=$(command -v rpm 2>/dev/null); HAVE_DPKG=$(command -v dpkg 2>/dev/null)
        for u in /etc/systemd/system/*.service /usr/lib/systemd/system/*.service; do
          [ -f "$u" ] || continue
          owned=1
          if [ -n "$HAVE_RPM" ]; then rpm -qf "$u" >/dev/null 2>&1 || owned=0
          elif [ -n "$HAVE_DPKG" ]; then dpkg -S "$u" >/dev/null 2>&1 || owned=0; fi
          # /etc/systemd/system legitimately holds admin units + symlinks
          if [ "$owned" = "0" ] && [ ! -L "$u" ]; then
            es=$(grep -E "^ExecStart=" "$u" 2>/dev/null | head -1)
            echo "  unowned unit (review): $u  | $es"
          fi
        done | head -25
        echo "  (admin-created units are legitimately unowned - judge by ExecStart)"

        echo ""
        echo "[*] === Lookalike / mimicry unit names ==="
        ls /etc/systemd/system/*.service /usr/lib/systemd/system/*.service 2>/dev/null | \
          xargs -n1 basename 2>/dev/null | \
          grep -iE "systemd-[a-z]*helper|dbus-?org|network-?manager-?[a-z]+|udev-?[a-z]+|kworker|kthread|kernel-?update|sys-?update" | \
          sort -u
        echo "  (compare against the real unit list; near-miss names are suspect)"

        echo ""
        echo "[*] === Recently created/modified units + timers ==="
        find /etc/systemd/system /usr/lib/systemd/system -name "*.service" -o -name "*.timer" 2>/dev/null | \
          xargs ls -lt 2>/dev/null | head -10

        echo ""
        echo "[*] === cron entries with odd payloads under blendy names ==="
        for c in /etc/cron.d/* /etc/cron.daily/* /etc/cron.hourly/* /var/spool/cron/* /var/spool/cron/crontabs/*; do
          [ -f "$c" ] || continue
          hit=$(grep -nE "/tmp/|/dev/shm/|curl|wget|base64|/var/tmp/|\| *(ba)?sh" "$c" 2>/dev/null)
          [ -n "$hit" ] && { echo "[FLAG] $c:"; echo "$hit" | head -3; }
        done

        echo ""
        echo "[*] === auditd: unit/cron writes outside package management ==="
        ausearch -k systemd_unit -i --start today 2>/dev/null | grep -ivE "comm=.(dpkg|rpm|yum|dnf|apt)" | tail -10
        ausearch -k cron_change -i --start today 2>/dev/null | grep -ivE "comm=.(dpkg|rpm)" | tail -10
    registry: |-
        Task / service masquerade artifacts:

        The idea:
          Persistence via systemd unit, timer, or cron is going to appear
          in an inventory. To survive review, the attacker NAMES the entry
          to look like a legitimate system service - systemd-update,
          dbus-org-helper, network-manager-dispatch, kernel-update.timer -
          so an analyst scanning systemctl or /etc/cron.d glides past it.

        The artifacts (same files as persistence, different lens):
          systemd units:
            /etc/systemd/system/*.service|*.timer   (admin/attacker units)
            /usr/lib/systemd/system/*               (package units - mimic targets)
            ~/.config/systemd/user/*                (user-scope persistence)
            /run/systemd/system/*                   (runtime/transient)
          cron:
            /etc/cron.d/* , /etc/cron.{daily,hourly,weekly}/*
            /var/spool/cron/* , /var/spool/cron/crontabs/*

        The detection principle (judge the target, not the name):
          A unit's NAME is attacker-controlled and meaningless for trust.
          Judge by:
            ExecStart PATH - a system-sounding unit whose ExecStart runs
              from /tmp, /dev/shm, /home, /var/tmp, or a hidden dir is
              malicious regardless of name
            PACKAGE OWNERSHIP - package-shipped units are owned by their
              package (rpm -qf/dpkg -S). An unowned, non-symlink unit in
              /usr/lib/systemd/system is suspect (admin units legitimately
              live in /etc/systemd/system, so weight that dir by ExecStart)
            BACKING BINARY - is the ExecStart binary itself package-owned
              and in a canonical path? (chains to T1036.005)

        High-signal indicators:
          a legit-sounding unit/timer with ExecStart in a writable path
          an unowned unit file in a package directory (/usr/lib/systemd)
          near-miss mimicry names (systemd-helperd, dbus-daemon-helper)
          Restart=always + ExecStart in /tmp (resilient malicious service)
          a cron entry under an innocuous name invoking curl|bash/base64
          recently-created unit on a host with no recent legitimate change

        Cross-reference:
          Persistence T1543.002 (systemd service) and T1053.003 (cron)
          cover the persistence mechanics; this row is the masquerading
          lens - the naming choice that hides those mechanics from review.
          The ExecStart binary is often itself name/location-masqueraded
          (T1036.005) and timestomped (T1070.006).
    ossdetect: |-
        Sigma rules:
        - file_event_lnx_systemd_unit_masquerade_name.yml
        - proc_creation_lnx_systemd_execstart_writable_path.yml
        - file_event_lnx_cron_masquerade_payload.yml

        Elastic detection rules:
        - Systemd Service with Suspicious ExecStart Path
        - Unowned Systemd Unit in Package Directory
        - Cron Job with Suspicious Command under Benign Name

        Auditd (same watches as persistence, read with the masquerade lens):
          -w /etc/systemd/system -p wa -k systemd_unit
          -w /usr/lib/systemd/system -p wa -k systemd_unit
          -w /etc/cron.d -p wa -k cron_change
          -w /var/spool/cron -p wa -k cron_change
          Exclude package managers; alert on shell-parented writes.

        Falco:
          rule: Create/modify systemd unit (non-package process)
          rule: Schedule cron job writing to suspicious path
          rule: Systemd service runs binary from writable directory

        The provenance check (scheduled hunt):
          for each unit: extract ExecStart binary; verify it is
          package-owned and in a canonical path; flag writable-path or
          unowned ExecStart. Ignore the unit NAME entirely.

        osquery:
          SELECT * FROM systemd_units;          (name, fragment_path)
          SELECT * FROM crontab;                (command, path)
          Join ExecStart/command targets to package ownership; baseline
          the expected unit set and diff.

        File integrity (FIM):
          AIDE/Tripwire/Wazuh on /etc/systemd/system, /usr/lib/systemd/
          system, /etc/cron.*, /var/spool/cron - alert on new entries;
          pair with an allow-list of expected units.

        Atomic Red Team:
          T1036.004 - masquerade task or service tests
    tools: |-
        Masquerading a task or service:

        Scheduled persistence is inevitably going to show up in an
        inventory - systemctl list-units, the contents of /etc/cron.d.
        So the attacker's evasion is not to hide the entry but to make
        it boring: a name that blends into the dozens of legitimate
        system services, so the reviewer's eye slides over it.

        Why naming-based hiding fails to the right check:
          The unit or cron NAME is fully attacker-controlled and carries
          zero trust value. The trustworthy signals are the TARGET and
          the PROVENANCE:
            - ExecStart / the cron command: where does it actually run?
              A "systemd-update.service" whose ExecStart is /tmp/.x is
              malicious no matter how plausible the name.
            - Package ownership: is the unit file shipped by a package?
              Is the ExecStart binary package-owned and in a canonical
              path? Unowned units in package directories, or ExecStart
              pointing at unowned binaries, are the real tells.
          So the hunt ignores the name entirely and evaluates the target.

        The judgement nuance:
          Admin-created systemd units legitimately live in
          /etc/systemd/system and are NOT package-owned - that is normal,
          not malicious. The discriminator there is the ExecStart path
          and binary provenance, not ownership of the unit file itself.
          In /usr/lib/systemd/system (the package directory), by contrast,
          an unowned non-symlink unit is genuinely anomalous.

        The full masquerade chain:
          A polished implant combines all three T1036 rows: a service
          named to blend (T1036.004), whose ExecStart binary is named
          like and placed like a system tool (T1036.005), with that
          binary's timestamp cloned from a neighbor (T1070.006). Pulling
          any one thread - ExecStart path, package ownership, or ctime
          mismatch - unravels the set.

        OT/ICS relevance:
          OT Linux services are typically a small, stable, well-known set.
          That makes a masqueraded extra service relatively easy to spot
          against a tight baseline - maintain an allow-list of expected
          units and alert on any addition, judged by ExecStart.

        Threat actor use:
          Naming systemd services and cron jobs to imitate system
          components is routine in Linux persistence - cryptomining crews
          (Kinsing, TeamTNT) and targeted operators alike create
          system-sounding units pointing at their payloads.
    actors:
        - APT41
        - Kimsuky
        - Kinsing
        - TeamTNT
    analyst_notes: 'Masquerading a task or service accepts that scheduled persistence will appear in an inventory and instead makes the entry boring: a systemd unit, timer, or cron job named to blend into the dozens of legitimate system services (systemd-update, dbus-org-helper, network-manager-dispatch, kernel-update.timer) so a reviewer scanning systemctl or /etc/cron.d glides past it. Naming-based hiding fails to the right check because the unit or cron name is fully attacker-controlled and carries zero trust value, while the trustworthy signals are the target and the provenance: the ExecStart or cron command (a systemd-update.service whose ExecStart is /tmp/.x is malicious no matter how plausible the name) and package ownership (whether the unit file is shipped by a package and whether the ExecStart binary is package-owned and in a canonical path). So the hunt ignores the name entirely and evaluates the target. An important judgement nuance is that admin-created systemd units legitimately live in /etc/systemd/system and are not package-owned, which is normal rather than malicious, so the discriminator there is the ExecStart path and binary provenance rather than ownership of the unit file itself, whereas in the package directory /usr/lib/systemd/system an unowned non-symlink unit is genuinely anomalous. A polished implant combines all three T1036 rows: a service named to blend (T1036.004) whose ExecStart binary is named and placed like a system tool (T1036.005) with that binary''s timestamp cloned from a neighbor (T1070.006), so pulling any one thread (ExecStart path, package ownership, or ctime mismatch) unravels the set. This is especially tractable in OT Linux, where services are typically a small, stable, well-known set, making a masqueraded extra service relatively easy to spot against a tight allow-list of expected units. The technique is routine in Linux persistence, used by cryptomining crews like Kinsing and TeamTNT and by targeted operators alike. This row is the masquerading lens on the same artifacts that persistence T1543.002 and T1053.003 cover mechanically.'
