# sigma/rules/host/defense_evasion/t1112_modify_registry_disable_defenses_fileless_storag_2.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/defense_evasion.js
title: 'T1112 - Modify Registry (disable defenses, fileless storage, hidden keys) (Fileless storage heuristic: a registry value '
id: 6b1d1a6b-646e-5c2f-aa38-2b4e4e6a03bf
status: experimental
description: 'Creating, modifying, or deleting registry values to weaken security (UAC EnableLUA=0, firewall, SmartScreen, Defender keys), to store payloads or configuration filelessly (encoded blobs under Run keys or custom keys to evade disk-based detection), or to hide artifacts using null-embedded key names that reg.exe and regedit cannot display; performed via reg.exe, PowerShell, reg save/restore, or direct API Query block: Fileless storage heuristic: a registry value set with an unusually long data field (base64/hex blob) under Run/RunOnce or app keys.'
references:
    - MITRE ATT&CK T1112
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.defense-evasion
    - attack.t1112
    - attack.g0016
    - attack.g0096
    - attack.g0032
    - attack.g0069
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        winlog.event_id: 13
        registry.data.strings|re: .{200,}
        registry.path:
            - '*"\CurrentVersion\Run"*'
            - '*"\Software\Classes\"*'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: high
tonk:
    fidelity: full
    source_file: host/js/data/defense_evasion.js
    domain: host
    technique: T1112
    technique_name: Modify Registry
    os: win
    query_block: 3 of 3
    original_query: |-
        winlog.event_id: 13 AND registry.data.strings: /.{200,}/
        AND registry.path: (*"\CurrentVersion\Run"* OR *"\Software\Classes\"*)
    sysmon: |-
        // Registry telemetry (Sysmon EID 12/13/14)
        //   12 - key create/delete
        //   13 - value set (the workhorse - covers most T1112)
        //   14 - key/value rename
        EventID=13 TargetObject ends with one of the sensitive values below.

        // Common defense-weakening value sets (Sysmon 13):
        //   *\Policies\System\EnableLUA = 0                (UAC off)
        //   *\Policies\System\ConsentPromptBehaviorAdmin = 0
        //   *\Windows Defender\DisableAntiSpyware = 1
        //   *\SmartScreenEnabled = Off
        //   *\Policies\...\DisableTaskMgr / DisableRegistryTools = 1
        //   *\System\CurrentControlSet\Services\<sec svc>\Start = 4

        // Fileless storage tells (Sysmon 13):
        //   a Run/RunOnce or custom value containing a long base64/hex blob
        //   values under HKCU\Software\Classes\... or odd app keys holding
        //   encoded script (read back by a loader - cross-ref persistence)

        // Tooling (Sysmon EID 1):
        Image=*\reg.exe CommandLine=(*add* OR *delete* OR *save* OR *restore* OR *import*)
        Image=(*\powershell.exe OR *\pwsh.exe) CommandLine=(*Set-ItemProperty*
          OR *New-ItemProperty* OR *Remove-ItemProperty* OR *reg save* OR *[Microsoft.Win32.Registry]*)
        Image=*\regini.exe
    powershell: |-
        # T1112 - Modify Registry hunt (defense-weakening + fileless + hidden)
        # Registry is BOTH a weapon (disable security) and a hiding place
        # (fileless payload/config, null-byte-hidden keys).

        Write-Host "[*] === Defense-weakening values ==="
        $checks = @(
          @{P='HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\Policies\System'; N='EnableLUA'; Bad=0},
          @{P='HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\Policies\System'; N='ConsentPromptBehaviorAdmin'; Bad=0},
          @{P='HKLM:\SOFTWARE\Policies\Microsoft\Windows Defender'; N='DisableAntiSpyware'; Bad=1},
          @{P='HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\Policies\Explorer'; N='DisableTaskMgr'; Bad=1}
        )
        foreach($c in $checks){
          $v = (Get-ItemProperty $c.P -Name $c.N -ErrorAction SilentlyContinue).$($c.N)
          if($v -ne $null -and $v -eq $c.Bad){ Write-Host ("[FLAG] {0}\{1} = {2} (weakened)" -f $c.P,$c.N,$v) }
        }

        Write-Host "[*] === SmartScreen state ==="
        (Get-ItemProperty 'HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\Explorer' -Name SmartScreenEnabled -ErrorAction SilentlyContinue).SmartScreenEnabled

        Write-Host "[*] === Image File Execution Options (debugger hijack / SilentProcessExit) ==="
        $ifeo='HKLM:\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Image File Execution Options'
        Get-ChildItem $ifeo -ErrorAction SilentlyContinue | ForEach-Object {
          $dbg=(Get-ItemProperty $_.PSPath -Name Debugger -ErrorAction SilentlyContinue).Debugger
          if($dbg){ Write-Host ("[FLAG] IFEO Debugger on {0} -> {1}" -f $_.PSChildName,$dbg) }
        }

        Write-Host "[*] === Fileless: oversized blobs under Run/RunOnce ==="
        foreach($rk in 'HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\Run','HKCU:\SOFTWARE\Microsoft\Windows\CurrentVersion\Run'){
          $p = Get-ItemProperty $rk -ErrorAction SilentlyContinue
          if($p){ $p.PSObject.Properties | Where-Object { $_.Value -is [string] -and $_.Value.Length -gt 200 } |
            ForEach-Object { Write-Host ("[FLAG] large value {0}\{1} ({2} chars)" -f $rk,$_.Name,$_.Value.Length) } }
        }

        Write-Host "[*] === Null-byte hidden keys (reg.exe/regedit cannot show these) ==="
        Write-Host "    reg.exe and Get-ItemProperty miss keys with embedded null"
        Write-Host "    bytes in the name. Use Autoruns (handles them) or the native"
        Write-Host "    RegEnumKeyEx API / SysinternalsSuite to enumerate honestly."

        Write-Host "[*] === Recent reg.exe / registry tooling from event log ==="
        Get-WinEvent -FilterHashtable @{LogName='Microsoft-Windows-Sysmon/Operational';Id=13} -MaxEvents 200 -ErrorAction SilentlyContinue |
          Where-Object { $_.Message -match 'EnableLUA|DisableAntiSpyware|SmartScreen|Image File Execution|DisableTaskMgr' } |
          Select-Object TimeCreated, @{n='Detail';e={($_.Message -split "`n" | Select-String 'TargetObject:')}} -First 20 | Format-Table -Auto
    registry: |-
        Modify Registry - the registry as weapon and hiding place:

        The registry is both a control surface (flip a value, weaken a
        defense) and a covert store (hold a payload or config with no file
        on disk). T1112 covers the EVASION uses; persistence Run-key abuse
        (T1547.001) and the Defender keys (T1562.001) are adjacent and
        often co-occur.

        Defense-weakening values (representative, not exhaustive):
          HKLM\...\Policies\System\EnableLUA = 0           (UAC disabled)
          HKLM\...\Policies\System\ConsentPromptBehaviorAdmin = 0
          HKLM\SOFTWARE\Policies\Microsoft\Windows Defender\DisableAntiSpyware=1
          HKLM\...\Explorer\SmartScreenEnabled = Off
          HKLM\...\Policies\Explorer\DisableTaskMgr / DisableRegistryTools = 1
          HKLM\SYSTEM\CurrentControlSet\Services\<security svc>\Start = 4
          Image File Execution Options\<exe>\Debugger   (hijack a binary's
            launch) and the GlobalFlag/SilentProcessExit evasion variants

        Fileless storage (evade disk-based detection):
          encoded payloads or config written as registry VALUES, read back
          and executed by a small loader. Classic homes: Run/RunOnce values
          holding base64, HKCU\Software\Classes\<random>, or bespoke app
          keys. The tell is an unusually long binary/string value where a
          short path or flag is expected - hunt value LENGTH and entropy.

        Hidden keys (anti-enumeration):
          keys/values created via the native API with an embedded NULL byte
          in the name are not displayed by reg.exe or regedit (which use the
          Win32 string APIs that stop at the null). Autoruns and direct
          RegEnumKeyEx-based tools DO show them. If a host behaves as though
          a Run entry exists but none is visible, suspect a null-byte key.

        Telemetry:
          Sysmon EID 13 (value set) is the workhorse, 12 (key create/delete)
          and 14 (rename) round it out. Baseline the sensitive keys above and
          alert on changes; reg save of SAM/SYSTEM/SECURITY hives is a
          credential-access co-signal (T1003) worth wiring in here too.

        Why this locks the Windows schema:
          this row exercises the Registry tab as ACTUAL registry (keys,
          values, hives, Sysmon 12/13/14) rather than the Linux File
          Artifacts repurposing - the structural difference between the
          Windows and Linux rows on this page.
    ossdetect: |-
        Sigma rules:
        - registry_set_disable_uac_enablelua.yml
        - registry_set_defender_disable.yml
        - registry_set_disable_smartscreen.yml
        - registry_set_ifeo_debugger.yml (Image File Execution Options)
        - registry_set_susp_run_key_blob.yml (fileless oversized value)

        Elastic detection rules:
        - UAC Disabled via Registry (EnableLUA)
        - Image File Execution Options Debugger Set
        - SmartScreen Disabled via Registry
        - Suspicious Large Registry Value (fileless storage)

        Sysmon config (the core telemetry):
          RegistryEvent 12/13/14 with a TargetObject allow/deny tuned to the
          sensitive keys (UAC, Defender, SmartScreen, services Start, IFEO,
          Run/RunOnce). Value-set (13) is the workhorse.

        Autoruns / Sysinternals:
          Autoruns enumerates persistence + IFEO keys AND handles null-byte
          hidden keys that reg.exe/regedit cannot display - run it for honest
          enumeration.

        osquery:
          SELECT path,name,data FROM registry WHERE path LIKE
            'HKEY_LOCAL_MACHINE\...\Policies\System\EnableLUA';
          SELECT * FROM registry WHERE path LIKE '%\Image File Execution Options\%'
            AND name='Debugger';
          schedule sensitive-key sweeps and diff against baseline.

        Baseline-and-diff:
          the sensitive value set is small and stable; periodic comparison
          against a known-good snapshot catches single-value weakening that
          signatures may miss.

        Atomic Red Team:
          T1112 - modify registry (UAC, Defender, IFEO, fileless) tests
    tools: |-
        Modifying the registry for evasion:

        On Windows the registry is where a lot of security posture lives, so
        it is simultaneously something an attacker breaks (to weaken
        defenses) and something they hide inside (to avoid the disk). T1112
        is the evasion lens on both.

        As a weapon - weakening defenses:
          a single value flip can disable UAC (EnableLUA=0), turn off
          SmartScreen, set a security service Start type to disabled, or
          enable a debugger hijack via Image File Execution Options. These
          are high-signal when you baseline the keys: there is a small,
          well-known set of "never legitimately changed at runtime" values,
          and a Sysmon EID 13 on any of them deserves an alert.

        As a hiding place - fileless storage:
          malware increasingly stores its payload or configuration as a
          registry value rather than a file, then uses a tiny loader to read
          and execute it. This defeats file-scanning and leaves a thinner
          forensic trail. The detection pivot is data SHAPE: a Run value or
          app key holding a long base64/hex blob where a short path belongs.
          Hunt on value length and entropy, not just on key path.

        The anti-enumeration trick - null-byte keys:
          because reg.exe and regedit use string APIs that terminate at a
          null, a key or value created with an embedded null in its name is
          invisible to them while still functioning. Autoruns and tools that
          call RegEnumKeyEx directly reveal them. This is the registry analogue
          of the Linux hidden-process trick: trust the structural enumerator,
          not the convenient tool that can be fooled.

        The detection backbone:
          Sysmon 12/13/14 plus a baselined watch-list of sensitive keys is
          the core. Layer Autoruns (for null-byte and persistence keys),
          osquery's registry table for scheduled sweeps, and an alert on reg
          save of the credential hives. Because much of this is one-value-
          deep, baseline-and-diff beats signature matching.

        Relationship to the rest of the chunk:
          T1112 overlaps the Defender row (those disables ARE registry
          writes) but extends to UAC/SmartScreen/IFEO weakening and to
          fileless storage and hidden keys - the broader registry-as-evasion
          surface. It is also the row that exercises the Windows Registry tab
          as real registry, setting the schema for every Windows row to come.

        Threat actor use:
          fileless registry storage is a hallmark of families like Kovter and
          Poweliks; UAC/Defender/SmartScreen weakening via registry is broad
          across ransomware and commodity malware; IFEO and SilentProcessExit
          abuse appears in more targeted tradecraft.
    actors:
        - APT29
        - APT41
        - Lazarus
        - MuddyWater
    analyst_notes: 'On Windows the registry is where a lot of security posture lives, so it is simultaneously something an attacker breaks to weaken defenses and something they hide inside to avoid the disk, and T1112 is the evasion lens on both. As a weapon a single value flip can disable UAC with EnableLUA=0, turn off SmartScreen, set a security service Start type to disabled, or enable a debugger hijack via Image File Execution Options, and these are high-signal when the keys are baselined because there is a small well-known set of values never legitimately changed at runtime, so a Sysmon EID 13 on any of them deserves an alert. As a hiding place, malware increasingly stores its payload or configuration as a registry value rather than a file and uses a tiny loader to read and execute it, which defeats file-scanning and leaves a thinner forensic trail, so the detection pivot is data shape: a Run value or app key holding a long base64 or hex blob where a short path belongs, hunted on value length and entropy rather than key path alone. The anti-enumeration trick is the null-byte key: because reg.exe and regedit use string APIs that terminate at a null, a key or value created with an embedded null in its name is invisible to them while still functioning, whereas Autoruns and tools that call RegEnumKeyEx directly reveal it, which is the registry analogue of the Linux hidden-process trick where you trust the structural enumerator rather than the convenient tool that can be fooled. The detection backbone is Sysmon 12/13/14 plus a baselined watch-list of sensitive keys, layered with Autoruns for null-byte and persistence keys, osquery''s registry table for scheduled sweeps, and an alert on reg save of the credential hives, and because much of this is one-value-deep, baseline-and-diff beats signature matching. T1112 overlaps the Defender row since those disables are registry writes but extends to UAC, SmartScreen, and IFEO weakening and to fileless storage and hidden keys, and it is the row that exercises the Windows Registry tab as real registry, setting the schema for every Windows row to come. Fileless registry storage is a hallmark of families like Kovter and Poweliks, registry-based weakening is broad across ransomware and commodity malware, and IFEO and SilentProcessExit abuse appears in more targeted tradecraft.'
