# sigma/rules/host/defense_evasion/t1562_001_disable_security_tools_and_edr_monitoring_agents_3.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/defense_evasion.js
title: T1562.001 - Disable Security Tools and EDR/Monitoring Agents (Package removal of a security product)
id: 0fc5faf7-2ef0-58c3-95ed-1771fd00d29e
status: experimental
description: 'Termination, stopping, or uninstallation of endpoint security and monitoring agents: killing or masking Falco, Wazuh, osquery, EDR sensors, ClamAV, or cloud-provider agents; cloud crypto crews specifically uninstall Alibaba/Tencent/AWS agents; or unloading kernel-based sensors (eBPF/kernel module) Query block: Package removal of a security product.'
references:
    - MITRE ATT&CK T1562.001
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.defense-evasion
    - attack.t1562.001
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        process.name:
            - rpm
            - apt
            - apt-get
            - dpkg
            - yum
            - dnf
        process.command_line|all:
            - '*-e*'
            - '*remove*'
            - '*purge*'
            - '*erase*'
            - '*falco*'
            - '*clamav*'
            - '*wazuh*'
            - '*osquery*'
            - '*falcon*'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/defense_evasion.js
    domain: host
    technique: T1562.001
    technique_name: 'Impair Defenses: Disable or Modify Tools'
    os: linux
    query_block: 4 of 5
    original_query: |-
        process.name: ("rpm" OR "apt" OR "apt-get" OR "dpkg" OR "yum" OR "dnf")
        AND process.command_line: ((*-e* OR *remove* OR *purge* OR *erase*)
          AND (*falco* OR *clamav* OR *wazuh* OR *osquery* OR *falcon*))
    sysmon: |-
        // Sysmon for Linux EID 1 - security agent disablement

        // Stopping / masking / disabling security services
        Image=(*/systemctl OR */service) CommandLine matches:
          *stop* *disable* *mask* combined with:
          *falco*  *wazuh*  *osquery*  *clamav*  *clamd*  *ossec*
          *crowdstrike* *falcon-sensor* *carbonblack* *cbagent*
          *sysmon* *auditd* *filebeat* *auditbeat* *sentinelone*

        // Killing agent processes
        Image=(*/kill OR */pkill OR */killall) CommandLine matches:
          *falco* *wazuh* *osqueryd* *clamd* *falcon* *cbagent*
          *aliyun* *AliYunDun* *aegis* *YunJing* *BcmAgent* (cloud agents)

        // Cloud-agent UNINSTALL (Rocke/TeamTNT signature behavior)
        CommandLine matches:
          */uninstall.sh*  *aegis*uninstall*  *AliYunDunUpdate*
          *quartz_uninstall*  *YunJing*  *bcm*uninstall*
          *rpm -e* / *apt remove* with a security package name

        // Unloading kernel/eBPF sensors
        Image=(*/rmmod OR */modprobe) CommandLine matches: *-r* (sensor module)

        // Auditd
        -w /etc/systemd/system -p wa -k svc_change
        -a always,exit -F arch=b64 -S kill -k kill_signal
        -a always,exit -F arch=b64 -S execve -F path=/usr/bin/pkill -k pkill_exec
    powershell: |-
        #!/bin/bash
        # T1562.001 - Security agent / EDR disablement hunt

        echo "[*] === Expected security agents: are they RUNNING? ==="
        for svc in falco wazuh-agent ossec osqueryd clamav-daemon clamd \
                   filebeat auditbeat falcon-sensor sysmon; do
          state=$(systemctl is-active "$svc" 2>/dev/null)
          enab=$(systemctl is-enabled "$svc" 2>/dev/null)
          [ "$state" = "unknown" ] && [ -z "$enab" ] && continue
          flag=""
          [ "$state" != "active" ] && flag="  <-- NOT ACTIVE"
          [ "$enab" = "masked" ] && flag="  <-- MASKED (blocked from starting)"
          echo "  $svc: active=$state enabled=$enab$flag"
        done

        echo ""
        echo "[*] === Cloud security agents (Rocke/TeamTNT uninstall targets) ==="
        for proc in aliyun AliYunDun aegis_cli AliHips YunJing BcmAgent \
                    qualys-cloud-agent amazon-ssm-agent; do
          if pgrep -f "$proc" >/dev/null 2>&1; then
            echo "  $proc: running"
          else
            # only flag if there's evidence it was once installed
            if ls /usr/local/aegis* /usr/local/qcloud* /opt/alibaba* 2>/dev/null | grep -q .; then
              echo "[FLAG] $proc: install dir present but process NOT running (uninstalled/killed?)"
            fi
          fi
        done

        echo ""
        echo "[*] === Masked units (blocked from ever starting) ==="
        systemctl list-unit-files --state=masked 2>/dev/null | grep -iE "falco|wazuh|clam|osquery|ossec|falcon|audit"

        echo ""
        echo "[*] === Recent service stop/disable in journal ==="
        journalctl --since today 2>/dev/null | grep -iE "Stopped|Stopping|Deactivated" | \
          grep -iE "falco|wazuh|clam|osquery|ossec|falcon|auditd" | tail -20

        echo ""
        echo "[*] === auditd: kills sent to agent PIDs + service edits ==="
        ausearch -k kill_signal -i --start today 2>/dev/null | \
          grep -iE "falco|wazuh|osquery|clam|falcon" | tail -15
        ausearch -k svc_change -i --start today 2>/dev/null | tail -10

        echo ""
        echo "[*] === Security package removal evidence ==="
        grep -hiE "remove|purge|erase" /var/log/dpkg.log /var/log/yum.log /var/log/dnf.log 2>/dev/null | \
          grep -iE "falco|clamav|wazuh|osquery|falcon|ossec" | tail -15
        echo "  (cross-check: agent install dir exists but package gone = uninstalled)"

        echo ""
        echo "[*] === eBPF/kernel sensor presence ==="
        lsmod 2>/dev/null | grep -iE "falco|sysdig|sensor" || echo "  no obvious sensor kmod"
        ls /sys/fs/bpf/ 2>/dev/null | head
    registry: |-
        Security agent / EDR disablement artifacts:

        Common Linux security agents (the targets):
          Detection/EDR: Falco, Wazuh, OSSEC, osquery/osqueryd,
            CrowdStrike falcon-sensor, SentinelOne, Carbon Black cbagent
          AV: ClamAV (clamd/clamav-daemon)
          Telemetry shippers: filebeat, auditbeat, fluentd
          Cloud-provider agents (heavily targeted by crypto crews):
            Alibaba: aegis / AliYunDun / AliHips (CloudMonitor, Security)
            Tencent: YunJing / BcmAgent (Yunjing host security)
            AWS: amazon-ssm-agent ; Qualys: qualys-cloud-agent

        Disablement methods:
          systemctl stop/disable/mask <agent>   (mask = cannot restart)
          kill / pkill / killall <agent process>
          rpm -e / apt remove / yum erase <agent package>
          cloud-agent uninstall scripts (the Rocke/TeamTNT signature):
            /usr/local/aegis/aegis_uninstall.sh
            AliYunDunUpdate, quartz_uninstall, YunJing uninstall
          rmmod / modprobe -r <kernel sensor module>
          config neutering (point agent at /dev/null, blank ruleset)

        The Rocke/TeamTNT cloud-agent TTP (well documented):
          These crews ship scripts that specifically detect and
          UNINSTALL Alibaba and Tencent cloud security agents before
          deploying miners. Presence of an agent's install directory
          WITHOUT its running process/package is the tell.

        High-signal indicators:
          an expected agent service in state inactive/masked
          agent install dir present but process absent (uninstalled)
          kill/pkill targeting a known agent process name
          security package in dpkg/yum remove logs
          masked security units (blocked from starting)
          agent stopped checking in to its console (heartbeat gap)
          cloud-agent uninstall script names in history/auditd

        The resilience principle:
          Agents that report to an external console make local
          disablement visible: the MISSING HEARTBEAT from a host that
          is otherwise alive is the detection. EDR consoles, Wazuh
          manager, and Falcon cloud all surface agent-down as an alert.
    ossdetect: |-
        Sigma rules:
        - lnx_security_tools_disabling.yml
        - proc_creation_lnx_clamav_service_stop.yml
        - proc_creation_lnx_cloud_agent_uninstall.yml
        - proc_creation_lnx_kill_security_process.yml

        Elastic detection rules:
        - Security Software Discovery then Termination
        - Cloud Security Agent Uninstallation (Alibaba/Tencent)
        - Tampering or Stopping of Falco / Wazuh / osquery

        Auditd:
          -w /etc/systemd/system -p wa -k svc_change
          -a always,exit -F arch=b64 -S kill -k kill_signal
          Correlate kill targets against known agent PIDs.

        Falco (catches kills of OTHER agents, and vice versa):
          rule: Terminate Security Tool / agent process
          rule: Disable Security Service
          Layered sensors mutually witness each other's disablement.

        THE control - external console / heartbeat monitoring:
          Wazuh manager: agent disconnected alert
          CrowdStrike/SentinelOne cloud: sensor-health / host-offline
          SIEM: missing-heartbeat rule per registered host
          Server-side agent-down alerting is unreachable by the host
          attacker and is the most reliable detection.

        osquery (inventory + drift):
          SELECT name FROM services WHERE name LIKE '%falco%' ...;
          Compare running services against a golden security baseline.

        Immutability / restart protection:
          systemd: Restart=always + protect units from masking via
          policy; file immutability (chattr +i) on agent configs.

        Atomic Red Team:
          T1562.001 - disable/modify tools (service stop, kill, uninstall)
    tools: |-
        Disabling security tools and agents:

        Before deploying a payload (miner, ransomware, backdoor), an
        adversary commonly clears the field by stopping or removing the
        host's security agents. On Linux servers this prominently
        includes CLOUD-PROVIDER agents, which is a defining behavior of
        several cryptomining crews.

        The signature cloud behavior (Rocke, TeamTNT, 8220):
          These actors ship uninstall routines that specifically target
          Alibaba Cloud (aegis/AliYunDun) and Tencent Cloud (YunJing)
          host-security agents. The scripts call the vendors' own
          uninstall paths to cleanly remove monitoring before mining.
          This is so consistent it is itself an attribution signal.

        General agent-killing:
          systemctl stop/mask, kill/pkill of falco/wazuh/osquery/clamd,
          package removal, or unloading kernel/eBPF sensors. Masking is
          favored over stopping because it prevents auto-restart.

        Why the cloud/console model wins:
          An on-host kill can blind the host locally, but agents that
          stream to an external manager (Wazuh manager, Falcon cloud,
          a SIEM) turn the kill into an alert: the agent-down / missing-
          heartbeat event fires from the SERVER side, which the attacker
          cannot reach. Local evasion, remote detection.

        The cross-tool principle:
          No single agent should be the only sensor. auditd + Falco +
          an EDR provide overlapping coverage; disabling one is caught
          by another (e.g. auditd records the pkill of Falco; Falco
          records the stop of auditd). This reference's hunts assume
          this layered model.

        Threat actor use:
          Rocke and TeamTNT cloud-agent uninstallation is extensively
          documented (Palo Alto Unit42, others). EDR/agent killing more
          broadly is near-universal pre-payload across crimeware and
          targeted Linux intrusions.
    actors:
        - Rocke
        - TeamTNT
        - 8220 Gang
    analyst_notes: 'Disabling security tools is a standard pre-payload step where an adversary clears the field before deploying a miner, ransomware, or backdoor, and on Linux servers it prominently includes cloud-provider agents, which is a defining behavior of several cryptomining crews. The signature cloud TTP, extensively documented for Rocke, TeamTNT, and 8220 Gang, is shipping uninstall routines that specifically target Alibaba Cloud (aegis/AliYunDun) and Tencent Cloud (YunJing) host-security agents by calling the vendors'' own uninstall paths, a behavior so consistent it serves as an attribution signal. More generally, adversaries stop, kill, mask, or remove agents like Falco, Wazuh, osquery, ClamAV, and commercial EDR sensors, with masking favored over stopping because it prevents auto-restart. The decisive defensive principle is the external-console model: an on-host kill blinds the host locally, but agents that stream to an external manager turn the kill into a server-side alert, because the agent-down or missing-heartbeat event fires from infrastructure the attacker cannot reach, so local evasion becomes remote detection. The complementary principle is layered sensors, where auditd, Falco, and an EDR provide overlapping coverage so that disabling one is witnessed by another, for example auditd recording the pkill of Falco while Falco records the stop of auditd. For hunting, the highest-signal indicators are an expected agent in inactive or masked state, an agent install directory present without its running process (the uninstall tell), kill or pkill targeting known agent process names, and security packages appearing in removal logs.'
