# sigma/rules/host/persistence/t1037_004_sysv_init_d_script_and_rc_d_symlink_persistence_2.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1037.004 - SysV init.d Script and rc.d Symlink Persistence (Registration commands)
id: 64adf5d8-42f8-52e5-9ae2-1d2634c89b31
status: experimental
description: 'A new or modified script in /etc/init.d/ with corresponding symlinks in /etc/rcN.d/ (or chkconfig/update-rc.d registration), particularly when the script is not owned by any installed package and its target binary resides in a writable or non-standard path Query block: Registration commands.'
references:
    - MITRE ATT&CK T1037.004
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1037.004
    - attack.g0096
    - attack.g0034
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        process.name:
            - update-rc.d
            - chkconfig
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: high
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1037.004
    technique_name: 'Boot or Logon Initialization Scripts: RC Scripts'
    os: linux
    query_block: 3 of 3
    original_query: 'process.name: ("update-rc.d" OR "chkconfig")'
    sysmon: |-
        // Auditd rules for init.d and rc.d writes
        -w /etc/init.d/ -p wa -k initd_write
        -w /etc/rc0.d/ -p wa -k rcd_write
        -w /etc/rc1.d/ -p wa -k rcd_write
        -w /etc/rc2.d/ -p wa -k rcd_write
        -w /etc/rc3.d/ -p wa -k rcd_write
        -w /etc/rc5.d/ -p wa -k rcd_write

        // update-rc.d or chkconfig invocation
        Image=*/update-rc.d OR Image=*/chkconfig
    powershell: |-
        # Hunt for init.d persistence
        echo "[*] === Non-package init.d scripts ==="
        for f in /etc/init.d/*; do
          [ -f "$f" ] || continue
          pkg=$(dpkg -S "$f" 2>/dev/null || rpm -qf "$f" 2>/dev/null)
          if [ -z "$pkg" ] || echo "$pkg" | grep -q 'not found|not owned'; then
            echo "  UNPACKAGED: $f"
            head -3 "$f"
            # Check what it runs
            grep -E '^[^#]*(/tmp/|/dev/shm/|/var/tmp/|curl|wget|python|perl)' "$f" 2>/dev/null
            echo "  ---"
          fi
        done

        echo ""
        echo "[*] === rc.d symlinks to non-package scripts ==="
        find /etc/rc?.d/ -type l 2>/dev/null | while read lnk; do
          target=$(readlink -f "$lnk")
          pkg=$(dpkg -S "$target" 2>/dev/null || rpm -qf "$target" 2>/dev/null)
          if [ -z "$pkg" ] || echo "$pkg" | grep -q 'not found|not owned'; then
            echo "  $lnk -> $target (UNPACKAGED)"
          fi
        done
    registry: |-
        No registry artifact (Linux technique).

        SysV init persistence mechanics:
        - /etc/init.d/<name> is the service script
        - /etc/rcN.d/S##<name> symlinks control start at runlevel N
        - /etc/rcN.d/K##<name> symlinks control stop at runlevel N
        - update-rc.d (Debian) or chkconfig (RHEL) manage the symlinks

        On systemd systems, SysV init scripts still work via the
        systemd-sysv-generator compatibility layer, which creates
        transient .service units for each init.d script at boot.
    ossdetect: |-
        Auditd:
        - -w /etc/init.d/ -p wa -k initd_write
        Wazuh:
        - FIM on /etc/init.d/ and /etc/rc*.d/
        osquery:
        - startup_items table
    tools: |-
        osquery:
          SELECT * FROM startup_items WHERE type='init';
        Wazuh FIM on /etc/init.d/ and /etc/rc*.d/
        Lynis: init script audit
    actors:
        - APT41
        - Sandworm
    analyst_notes: 'SysV init.d persistence is a legacy mechanism that remains effective because (a) systemd''s sysv-generator compatibility layer still runs these scripts on modern systems, and (b) many defenders only audit systemd units, not init.d. The same package-verification approach works: enumerate all init.d scripts, verify against the package manager, flag unpackaged scripts. Scripts pointing at binaries in /tmp, /dev/shm, or user home directories are high-confidence malicious. On ESXi hosts, /etc/rc.local.d/ is the primary persistence mechanism since ESXi uses a custom init, not systemd.'
