# sigma/rules/host/persistence/t1053_005_hidden_scheduled_tasks_sd_manipulation_and_tree_2.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1053.005 - Hidden Scheduled Tasks (SD Manipulation and Tree Hiding) (Security Event 4698 created)
id: e04cb6ed-3880-5bad-8f43-4ffffb2e1f8d
status: experimental
description: 'Tasks present in C:\Windows\System32\Tasks XML but missing from schtasks.exe enumeration - registry/filesystem diff reveals hidden tasks Query block: Security Event 4698 created.'
references:
    - MITRE ATT&CK T1053.005
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1053.005
    - attack.g0125
    - attack.g0096
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 4698
        winlog.channel: Security
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1053.005
    technique_name: 'Scheduled Task/Job: Scheduled Task (Persistence Angle)'
    os: win
    query_block: 3 of 3
    original_query: |-
        winlog.event_id: 4698
        AND winlog.channel: "Security"
    sysmon: |-
        // Hidden task creation typically involves direct
        // registry/filesystem manipulation rather than
        // schtasks.exe, so process-based detection misses it.

        // EID 11 (FileCreate) catching task XML write:
        EventID=11
        TargetFilename=C:\Windows\System32\Tasks\*

        // EID 13 catching TaskCache registry manipulation:
        EventID=13
        TargetObject matches:
          *\Schedule\TaskCache\Tree\*
          OR *\Schedule\TaskCache\Tasks\*

        // Cross-correlate: XML written without corresponding
        // schtasks.exe process creation = direct manipulation

        // Security Event 4698 (created) WITHOUT a paired
        // schtasks.exe Sysmon EID 1 is unusual:
        EventID=4698
        LogName=Security
        // Then look for absence of corresponding EID 1
    powershell: |-
        # Hunt for hidden scheduled tasks - filesystem vs SCM diff

        # Method 1: Tasks in filesystem but not visible to Get-ScheduledTask
        $visibleTasks = Get-ScheduledTask | ForEach-Object {
          ($_.TaskPath + $_.TaskName).TrimStart('\\')
        }

        $fsTasks = Get-ChildItem 'C:\Windows\System32\Tasks' -Recurse -File |
          ForEach-Object {
            $_.FullName.Replace('C:\Windows\System32\Tasks\','').Replace('\\','\\')
          }

        $hidden = $fsTasks | Where-Object { $_ -notin $visibleTasks }

        if ($hidden) {
          Write-Host "Potentially hidden tasks (filesystem but not enumerable):" -ForegroundColor Red
          $hidden | ForEach-Object {
            $path = "C:\Windows\System32\Tasks\$_"
            if (Test-Path $path) {
              [PSCustomObject]@{
                TaskPath = $_
                XmlFile = $path
                Created = (Get-Item $path).CreationTime
                Modified = (Get-Item $path).LastWriteTime
                Size = (Get-Item $path).Length
              }
            }
          }
        } else {
          Write-Host "No hidden tasks detected via filesystem/SCM diff" -ForegroundColor Green
        }

        # Method 2: Check TaskCache registry for orphaned entries
        $registryGuids = Get-ChildItem 'HKLM:\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\Tasks' |
          Select-Object -ExpandProperty PSChildName

        $treeGuids = @()
        function Get-AllTreeTasks($key) {
          Get-ChildItem $key -EA SilentlyContinue | ForEach-Object {
            $id = (Get-ItemProperty $_.PSPath -Name Id -EA SilentlyContinue).Id
            if ($id) { $script:treeGuids += $id }
            Get-AllTreeTasks $_.PSPath
          }
        }
        Get-AllTreeTasks 'HKLM:\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\Tree'

        # GUIDs in Tasks key but not in Tree = orphaned/hidden tasks
        $orphans = $registryGuids | Where-Object { $_ -notin $treeGuids }
        if ($orphans) {
          Write-Host "Orphaned task GUIDs (in Tasks key but not in Tree):" -ForegroundColor Red
          $orphans
        }

        # Method 3: Read raw task XML for tasks that Get-ScheduledTask
        # may not fully surface
        Get-ChildItem 'C:\Windows\System32\Tasks' -Recurse -File | ForEach-Object {
          try {
            $xml = [xml](Get-Content $_.FullName -Raw)
            [PSCustomObject]@{
              Path = $_.FullName
              Author = $xml.Task.RegistrationInfo.Author
              RunAs = $xml.Task.Principals.Principal.UserId
              Trigger = ($xml.Task.Triggers.ChildNodes.Name) -join ', '
              Action = ($xml.Task.Actions.Exec.Command + ' ' + $xml.Task.Actions.Exec.Arguments)
            }
          } catch { }
        }
    registry: |-
        Scheduled task storage structure (Windows 10+):

        Filesystem:
        C:\Windows\System32\Tasks\<TaskPath>\<TaskName>
        - XML file with full task definition
        - This is the canonical task storage

        Registry mirror:
        HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Schedule\
          TaskCache\
            Tree\<TaskPath>\<TaskName>\
              Id = {GUID}  (links to Tasks key below)
              Index = <byte>
            Tasks\{GUID}\
              Actions = REG_BINARY  (action data)
              Triggers = REG_BINARY (trigger data)
              DynamicInfo = REG_BINARY (last run/result)
              Path = <task path back-reference>

        Hiding techniques:

        1. Filesystem only, no registry entry:
           - Write task XML to C:\Windows\System32\Tasks
           - Skip the registry TaskCache entries
           - Get-ScheduledTask uses registry - misses it
           - Task may or may not run depending on Windows version
             (newer Windows requires both)

        2. Registry only, no filesystem XML:
           - Write to TaskCache\Tasks\{GUID}
           - Reverses the above hiding pattern
           - Less common but documented in some malware

        3. Security descriptor manipulation:
           - HKLM\...\Schedule\TaskCache\Tree\<task>\
             SD = REG_BINARY (security descriptor)
           - Modify SD to deny enumeration access
           - Get-ScheduledTask and schtasks /query miss it
           - Task still runs because Task Scheduler service
             uses a different privilege level

        4. Tree entry hidden, GUID entry present:
           - Delete the Tree\<TaskName>\Id value
           - Task data still in Tasks\{GUID} - task runs
           - Tree enumeration misses the entry
           - PowerShell COM API and schtasks /query both miss it

        Detection approaches:

        - Diff filesystem vs Get-ScheduledTask output
          (shown in script above)
        - Check TaskCache\Tasks GUIDs against Tree entries
          (shown in script above)
        - Inspect task SD values for non-default ACLs:
          Get-Item 'HKLM:\...\TaskCache\Tree\<name>' |
          Get-Acl | Select-Object -Expand AccessToString

        The XML-based hunt approach reads C:\Windows\System32\
        Tasks directly and bypasses any registry-based hiding -
        even tasks hidden via SD manipulation in the registry
        remain visible if their XML still exists.
    ossdetect: |-
        Sigma:
        - file_event_win_scheduled_task_xml_creation.yml
        - registry_event_taskcache_modification.yml
        - registry_event_taskcache_sd_modification.yml

        Atomic Red Team:
        - T1053.005 has tests for task creation variants
        - Hidden task variant less commonly automated

        Hayabusa:
        - ScheduledTaskXmlWrite rules (EID 11)
        - TaskCacheModification rules (EID 13)

        Velociraptor:
        - Windows.System.ScheduledTasks
          (parses XML directly - surfaces hidden tasks
          that registry-based tools miss)
        - Best single tool for fleet-wide hunting

        Microsoft references:
        - 'Tarrask malware uses scheduled tasks for defense
          evasion' (Microsoft Threat Intelligence, 2022)
        - Foundational reference for SD-based task hiding
        - HAFNIUM tradecraft documentation

        KAPE (Eric Zimmerman):
        - Targets/Compound/ScheduledTasks.tkape
          (collects both XML files and registry hive)
        - Enables offline forensic analysis of
          both visible and hidden tasks
    tools: |-
        Hidden scheduled task tradecraft:

        Known adversary toolkits:
        - Various custom implants used by capable operators
        - Reported in Mandiant / CrowdStrike intrusion analyses
        - Documented use by Tarrask malware family
          (Microsoft 2022 disclosure on HAFNIUM tradecraft)

        Tarrask specifically (HAFNIUM):
        - Created scheduled tasks with deleted SD security
          descriptor in registry
        - Tasks invisible to schtasks.exe and services.msc
        - Persisted for months in long-dwell intrusions
        - Microsoft's 2022 disclosure brought the technique
          into mainstream defender awareness

        Standard implementation approach:
        1. Create task via schtasks.exe or COM API (visible)
        2. Modify TaskCache\Tree\<name>\SD value to deny
           enumeration access OR delete the SD value entirely
        3. Optional: also remove the Tree\<name>\Id mapping
           so registry-based enumeration misses it
        4. Task XML remains in filesystem - the run mechanism
           uses it - but standard enumeration tools don't
           surface the task

        Detection note:
        - This is sophisticated tradecraft - typically
          associated with APT-level operators
        - Commodity malware rarely bothers with this
        - Finding hidden tasks is high-confidence signal
          of capable / deliberate compromise

        Defensive recommendation:
        - Always cross-reference filesystem (System32\Tasks)
          against schtasks.exe / Get-ScheduledTask output
        - Use Velociraptor or similar XML-based enumeration
          rather than relying on the registry path alone
    actors:
        - HAFNIUM
        - APT41
    analyst_notes: 'Hidden scheduled tasks via SD manipulation came to mainstream attention in April 2022 when Microsoft disclosed the Tarrask malware family used by HAFNIUM (China-nexus). The technique exploits the gap between how tasks are stored (filesystem XML + registry TaskCache mirror) and how they''re typically enumerated (schtasks.exe and Get-ScheduledTask use the registry). By modifying or deleting the security descriptor on the TaskCache\Tree registry entry, the adversary makes the task invisible to standard tools while leaving the XML in place so it still runs. The filesystem-vs-SCM diff in the PowerShell script is the most reliable detection because it bypasses the registry entirely - if the XML exists in C:\Windows\System32\Tasks but Get-ScheduledTask doesn''t list it, something is hiding it. This is one of those techniques that''s rare enough that most defenders don''t actively look for it, which makes it disproportionately valuable for adversaries who do use it. Tarrask persisted in HAFNIUM intrusions for months in some cases before being detected. Pair this hunt with the standard task enumeration in the previous indicator: together they cover both visible-but-suspicious tasks and tasks that have been deliberately hidden. Velociraptor''s Windows.System.ScheduledTasks is the best fleet-wide tool because it reads XML directly rather than depending on the registry path.'
