# sigma/rules/host/persistence/t1053_003_cron_baseline_sweep_and_anomaly_hunt_hash_based.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: 'T1053.003 - Cron Baseline Sweep and Anomaly Hunt (Hash-based: alert when known cron file hash changes)'
id: 0e135cd6-99f4-55da-9a49-3edb7cdf2037
status: experimental
description: 'Enumerate all cron locations fleet-wide and diff against known-good baseline; alert on new @reboot, every-minute, or suspicious-path entries added outside change windows Query block: Hash-based: alert when known cron file hash changes.'
references:
    - MITRE ATT&CK T1053.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1053.003
logsource:
    category: file_event
    product: windows
detection:
    selection:
        event.module: file_integrity
        file.path: /etc/crontab
    filter_0:
        file.hash.sha256: <your-baseline-hash>
    condition: selection and not filter_0
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1053.003
    technique_name: 'Scheduled Task/Job: Cron'
    os: linux
    query_block: 3 of 4
    original_query: |-
        event.module: "file_integrity"
        AND file.path: "/etc/crontab"
        AND NOT file.hash.sha256: "<your-baseline-hash>"
    sysmon: |-
        # Generate and compare cron baseline
        # Run on all Linux hosts via orchestration tool

        # --- Snapshot all cron content ---
        # Output format: hostname | location | owner | content

        snapshot_cron() {
          local host=$(hostname)
          for f in /etc/crontab /etc/cron.d/* 
                    /etc/cron.hourly/* /etc/cron.daily/* 
                    /etc/cron.weekly/* /etc/cron.monthly/* 
                    /var/spool/cron/crontabs/* /var/spool/cron/*; do
            [ -f "$f" ] || continue
            stat -c "%h|%U|%G|%a|%Y" "$f" | 
              awk -v h="$host" -v p="$f" 
              '{print h"|"p"|"$0}'
          done
        }

        # --- Auditd real-time for @reboot hunting ---
        # Look for @reboot writes that were not from package manager:
        ausearch -k cron_persist --start today 2>/dev/null | 
          grep -A 5 "SYSCALL" | 
          grep -B 3 "@reboot"

        # --- Cross-host comparison (if central logging) ---
        # Collect cron file hashes for baseline comparison:
        find /etc/cron* /var/spool/cron/ -type f 2>/dev/null | 
          sort | xargs sha256sum 2>/dev/null
    powershell: |-
        #!/bin/bash
        # T1053.003 - Baseline and anomaly detection
        # Generates sha256 baseline of all cron content
        # Compare against previous run to detect changes

        BASELINE="/var/lib/security/cron_baseline.sha256"
        CURRENT="/tmp/cron_current_$(date +%s).sha256"
        ALERT_LOG="/var/log/security/cron_changes.log"

        mkdir -p "$(dirname "$BASELINE")" "$(dirname "$ALERT_LOG")"

        # Build current state hash
        find /etc/cron* /var/spool/cron/ -type f 2>/dev/null | sort | 
          while read f; do
            sha256sum "$f" 2>/dev/null
            echo "--- content of $f ---"
            cat "$f" 2>/dev/null
            echo ""
          done | sha256sum > "$CURRENT"

        # If baseline exists, compare
        if [ -f "$BASELINE" ]; then
          if ! diff -q "$CURRENT" "$BASELINE" > /dev/null 2>&1; then
            echo "[ALERT] $(date): Cron content changed" >> "$ALERT_LOG"
            # Deep comparison
            find /etc/cron* /var/spool/cron/ -type f 2>/dev/null | sort | 
              while read f; do
                sha256sum "$f"
              done > /tmp/cron_files_now.txt
            diff <(cat "$BASELINE") /tmp/cron_files_now.txt >> "$ALERT_LOG"
            echo "[*] Changed files:"
            diff <(cat "$BASELINE") /tmp/cron_files_now.txt

            # Specifically call out @reboot and downloader patterns
            echo "[*] @reboot entries:"
            grep -rh "@reboot" /etc/cron* /var/spool/cron/ 2>/dev/null
            echo "[*] Suspicious patterns:"
            grep -rEh "(curl|wget|bash|python|/tmp/|/dev/shm/|base64)" 
              /etc/cron* /var/spool/cron/ 2>/dev/null
          else
            echo "[OK] Cron baseline unchanged"
          fi
        else
          echo "[*] No baseline found - creating new baseline"
        fi

        cp "$CURRENT" "$BASELINE"
        rm -f "$CURRENT"
    registry: |-
        Baseline comparison targets:

        All paths in previous indicator plus:

        Access control (check for unexpected removal):
          /etc/cron.allow   - if present, limits who can use cron
          /etc/cron.deny    - denylist; adversaries may DELETE this file
                              to allow a previously denied account to cron

        Cron binary itself (watch for replacement):
          /usr/sbin/crond   - RHEL daemon binary
          /usr/sbin/cron    - Debian daemon binary
          Check hash against package manager:
            rpm -V cronie  (RHEL)
            dpkg --verify cron  (Debian)

        Log file for audit review:
          /var/log/syslog    - cron execution log on Debian
          /var/log/cron      - cron execution log on RHEL
          /var/log/cron.log  - some distros

        Integrity baseline approach:
          AIDE or Tripwire configured to watch:
            /etc/cron.d r+sha256+inode
            /etc/crontab r+sha256+inode
            /var/spool/cron/crontabs r+sha256+inode
    ossdetect: |-
        Velociraptor:
        - Linux.Sys.Crontab (primary fleet-sweep artifact)
        - Returns all cron entries with structured output
        - Can diff against uploaded baseline JSON

        osquery:
        - SELECT * FROM crontab;
        - SELECT * FROM file WHERE directory IN
          ('/etc/cron.d','/etc/cron.hourly',...)
          AND mtime > (strftime('%s','now') - 86400);

        AIDE (filesystem integrity):
        - Configure to watch /etc/cron*, /var/spool/cron
        - aide --check to compare against baseline db
        - aide --init to create initial baseline

        Tripwire:
        - Add cron directories to tripwire.cfg
        - Runs nightly; reports on any changes

        Auditd + ausearch:
        - ausearch -k cron_persist --start recent
        - aureport --file -ts today (file access summary)

        Sigma:
        - file_event_lnx_crontab_file_modification.yml
          (catches all listed cron paths)
    tools: |-
        Baseline sweep approach is critical when:
        - Real-time auditd was not deployed at compromise time
        - Attacker established persistence BEFORE logging was enabled
        - Investigating host with unknown dwell time

        Fleet-wide sweep tools:
        - Velociraptor: Linux.Sys.Crontab VQL artifact
          (enumerate, collect, compare fleet-wide)
        - osquery: crontab table (SELECT * FROM crontab)
        - Ansible/Salt: run cron snapshot playbook fleet-wide
        - PANIX --cron: simulates all cron persistence vectors
          for testing your baseline detection

        Known malware cron entry patterns:
        - Rocke: /etc/cron.d/apache  (spoofed name, curl | bash)
        - TeamTNT: every-minute entries downloading from 
          cloud-provider-spoofed C2 domains
        - Generic miner: /tmp/ or /var/tmp/ payload paths
        - ESXiArgs: /var/spool/cron/crontabs/root directly

        Hunting additional @reboot-equivalent persistence:
          systemctl list-timers --all
          (systemd timers: see T1053.006 for full coverage)
    actors:
        - Rocke
        - TeamTNT
    analyst_notes: 'The sweep approach is essential when you cannot trust that real-time detection caught every persistence write - either because auditd wasn''t running or because the adversary disabled logging temporarily. Key anomalies to prioritize: (1) @reboot entries not associated with a known package installation, (2) entries with paths in /tmp, /dev/shm, /var/tmp, or user home directories, (3) entries with download commands (curl, wget, python -c), (4) frequency higher than daily for root-owned entries (legitimate system jobs almost never need every-minute execution), (5) cron entries that reference files with names mimicking system components (systemd-update, network-manager, atd). The cron.allow / cron.deny pair is frequently overlooked: if cron.deny is deleted and cron.allow does not exist, ALL users can create crontabs. Adversaries who gain a low-privilege shell may delete cron.deny to allow their compromised user account to schedule tasks. Baseline the existence and content of both files.'
