# sigma/rules/host/persistence/t1053_005_tasks_with_persistence_intent_triggers_short_int.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1053.005 - Tasks with Persistence-Intent Triggers (Short-interval recurring tasks (beacon persistence))
id: 9e33aff0-5f90-5a24-8802-8524de13a6ca
status: experimental
description: 'Enumerate scheduled tasks and filter to those with OnLogon, AtStartup, OnEvent, or short-interval triggers - the persistence-specific subset Query block: Short-interval recurring tasks (beacon persistence).'
references:
    - MITRE ATT&CK T1053.005
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1053.005
    - attack.g0016
    - attack.g0096
    - attack.g0032
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 4698
        winlog.event_data.TaskContent:
            - '*PT1M*'
            - '*PT5M*'
            - '*PT10M*'
            - '*PT15M*'
    condition: selection
falsepositives:
    - The combination 'OnLogon trigger + action invokes powershell + action contains EncodedCommand or AppData path' is high-confidence malicious - any one of those alone has false positives, but the combination is rarely legitimate
level: high
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1053.005
    technique_name: 'Scheduled Task/Job: Scheduled Task (Persistence Angle)'
    os: win
    query_block: 3 of 3
    original_query: |-
        winlog.event_id: 4698
        AND winlog.event_data.TaskContent: (*PT1M* OR *PT5M* OR *PT10M* OR *PT15M*)
    sysmon: |-
        // This indicator is hunting-oriented - PowerShell
        // enumeration of established tasks rather than
        // real-time creation detection (T1053.005 in Execution).
        //
        // Sysmon contribution: historical task creation events:
        EventID=1
        Image=*\schtasks.exe
        CommandLine=*/create*
        // Cross-reference with /sc ONLOGON, /sc ONSTART,
        // /sc ONEVENT to find persistence-intent creations

        // Security Event 4698 (task created) - broader coverage:
        EventID=4698
        LogName=Security
        // TaskContent XML field reveals trigger type
        // Look for: <LogonTrigger>, <BootTrigger>,
        // <EventTrigger>, <TimeTrigger Repetition Interval="PT<short>">
    powershell: |-
        # Hunt scheduled tasks with persistence-intent triggers
        # OnLogon, AtStartup, OnEvent, or short recurring interval

        $persistenceTasks = @()

        Get-ScheduledTask | ForEach-Object {
          $task = $_
          $triggerTypes = @()
          $suspicious = @()

          foreach ($trigger in $task.Triggers) {
            $type = $trigger.CimClass.CimClassName

            # Classify the trigger
            switch ($type) {
              'MSFT_TaskLogonTrigger'   { $triggerTypes += 'OnLogon' }
              'MSFT_TaskBootTrigger'    { $triggerTypes += 'AtStartup' }
              'MSFT_TaskEventTrigger'   { $triggerTypes += 'OnEvent' }
              'MSFT_TaskSessionStateChangeTrigger' { $triggerTypes += 'SessionStateChange' }
              'MSFT_TaskTimeTrigger'    { $triggerTypes += 'Time' }
              'MSFT_TaskDailyTrigger'   { $triggerTypes += 'Daily' }
              'MSFT_TaskRegistrationTrigger' { $triggerTypes += 'OnRegistration' }
            }

            # Check for short recurring intervals (beacon-like)
            if ($trigger.Repetition.Interval) {
              $interval = $trigger.Repetition.Interval
              if ($interval -match 'PT([0-9]+)M' -and [int]$matches[1] -le 15) {
                $suspicious += "Short repetition interval ($interval)"
              }
              if ($interval -match 'PT([0-9]+)S' -and [int]$matches[1] -le 600) {
                $suspicious += "Sub-minute repetition ($interval)"
              }
            }
          }

          # Persistence-intent triggers (any of these alone qualifies)
          $persistenceTriggers = @('OnLogon','AtStartup','OnEvent','SessionStateChange')
          $isPersistence = $false
          foreach ($t in $triggerTypes) {
            if ($persistenceTriggers -contains $t) { $isPersistence = $true; break }
          }

          if (-not $isPersistence -and $suspicious.Count -eq 0) { return }

          # Pull the action details
          $actions = $task.Actions | ForEach-Object {
            if ($_.Execute) { "$($_.Execute) $($_.Arguments)" }
          }

          # Check action for additional risk signals
          $actionStr = $actions -join ' | '
          if ($actionStr -match '(powershell|cmd\.exe|wscript|cscript|mshta|rundll32|regsvr32)') {
            $suspicious += "Action invokes interpreter"
          }
          if ($actionStr -match '(AppData|Temp|ProgramData|\\Users\\Public)') {
            $suspicious += "Action in user-writable path"
          }

          $persistenceTasks += [PSCustomObject]@{
            TaskName    = $task.TaskName
            TaskPath    = $task.TaskPath
            Triggers    = ($triggerTypes -join ', ')
            Action      = $actionStr
            RunAs       = $task.Principal.UserId
            State       = $task.State
            Author      = $task.Author
            LastRunTime = $task.LastRunTime
            Suspicious  = ($suspicious -join '; ')
          }
        }

        $persistenceTasks | Sort-Object TaskPath | Format-Table -AutoSize
    registry: |-
        Scheduled task storage on disk:

        C:\Windows\System32\Tasks\<TaskPath>\<TaskName>
        - XML file containing full task definition
        - Human-readable - inspect with Get-Content

        Key XML elements for persistence analysis:
        <Triggers>
          <LogonTrigger>     <!-- persistence intent -->
            <UserId>...</UserId>
          </LogonTrigger>
          <BootTrigger>      <!-- persistence intent -->
          </BootTrigger>
          <EventTrigger>     <!-- can be persistence -->
            <Subscription>...event XPath query...</Subscription>
          </EventTrigger>
          <CalendarTrigger>  <!-- scheduled, may not be persistence -->
            <Repetition>
              <Interval>PT5M</Interval>  <!-- short = suspicious -->
            </Repetition>
          </CalendarTrigger>
        </Triggers>

        <Principal>
          <UserId>S-1-5-18</UserId>  <!-- SYSTEM = high privilege -->
          <RunLevel>HighestAvailable</RunLevel>
        </Principal>

        Registry equivalent locations:
        HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\
          Schedule\TaskCache\Tree\<TaskPath>\<TaskName>
        HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\
          Schedule\TaskCache\Tasks\{GUID}\
        - TaskCache\Tree maps task names to GUIDs
        - TaskCache\Tasks contains the action data per GUID
        - Useful for finding tasks that exist in registry
          but not in C:\Windows\System32\Tasks (rare,
          indicates manipulation)

        Cross-reference with Execution tactic:
        - T1053.005 in Execution covers task CREATION
          (real-time schtasks.exe, Security EID 4698)
        - This indicator covers HUNTING established tasks
        - Use both for full coverage

        Last-modified timestamps:
        (Get-Item "C:\Windows\System32\Tasks\<name>").LastWriteTime
        - Useful for incident timelining
        - Tells you when the task was installed or modified
    ossdetect: |-
        Sigma:
        - Real-time detection rules cover creation
        - Hunting established tasks is a manual workflow

        Atomic Red Team:
        - T1053.005 (multiple persistence trigger variants)
        - T1053.005 Test #1, #6 (OnLogon and OnStartup tasks)
        - Use to validate the hunt: install via Atomic,
          run PowerShell sweep, confirm detection

        Hayabusa:
        - ScheduledTaskCreation rules (real-time, EID 4698)
        - For hunting, use Get-ScheduledTask or schtasks /query

        Velociraptor:
        - Windows.System.ScheduledTasks
          (best fleet-wide artifact - enumerates all tasks
          with trigger and action detail)
        - Windows.Forensics.Autoruns

        Sysinternals autoruns.exe:
        - Scheduled Tasks tab
        - Highlights non-Microsoft signed tasks
        - VirusTotal lookup integration
        - For manual IR triage

        PowerShell native (the script above):
        - No tool deploy required
        - Works on any Windows with PowerShell 3.0+
        - Returns parsed trigger type and timing info
    tools: |-
        Persistence-intent trigger types:

        OnLogon (most common):
        - Fires when specified user logs in
        - "Any user" or specific account
        - Most common adversary persistence trigger
        - Standard Cobalt Strike persistence default

        AtStartup (system boot):
        - Fires when system boots
        - Requires SYSTEM-level task (admin to create)
        - Stronger persistence: runs before user logon

        OnEvent (event-driven):
        - Fires when specific Windows event log entry appears
        - Adversaries use specific event IDs as triggers
        - Example: trigger on Event ID 4624 LogonType=2
          (interactive logon - run code when someone logs in)
        - More evasive than OnLogon - less commonly monitored

        SessionStateChange:
        - Fires on logon/logoff/lock/unlock
        - Less common but documented

        Short recurring intervals:
        - PT1M, PT5M, PT10M repetitions
        - Used for beacon-style persistence
        - Adversary runs a check-in script every N minutes
        - Combined with stub task that contacts C2

        Adversary frameworks:
        SharPersist (/t schtask)
        Cobalt Strike (built-in task persistence)
        Empire / Starkiller modules
        Metasploit persistence
        Most commodity malware uses simpler schtasks /create
          with /sc ONLOGON

        Task disguise patterns (covered in T1053.005 Execution):
        - Names mimicking Windows tasks: "Windows Update",
          "Microsoft Office Health", "Google Update Task"
        - Task paths in \Microsoft\Windows\ to blend with
          legitimate Windows tasks
        - Author field set to "Microsoft Corporation"
          (easily spoofed, low fidelity by itself)
    actors:
        - APT29
        - APT41
        - Lazarus
    analyst_notes: 'This is the persistence-hunting counterpart to T1053.005 in Execution. The Execution side covers task CREATION (real-time schtasks.exe and Security EID 4698); this side covers HUNTING established tasks with persistence-specific characteristics. The detection narrows by trigger type: OnLogon, AtStartup, OnEvent, and SessionStateChange are the persistence-intent trigger families. Time-based triggers can be persistence (daily/weekly recurring) but are also legitimate for many scheduled jobs - they need additional context (action analysis) to classify. Short recurring intervals (PT1M-PT15M) are a strong signal for beacon-style persistence: legitimate tasks rarely repeat that frequently, while adversary beacons commonly do. The combination ''OnLogon trigger + action invokes powershell + action contains EncodedCommand or AppData path'' is high-confidence malicious - any one of those alone has false positives, but the combination is rarely legitimate. Pair this hunt with the Execution side''s third indicator (task action hunting) for full coverage: this indicator filters by trigger type, that one filters by action path - together they catch the persistence-intent + suspicious-execution combination from both angles.'
