# sigma/rules/host/persistence/t1543_003_hidden_services_via_security_descriptor_manipula_1.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1543.003 - Hidden Services via Security Descriptor Manipulation (sc.exe used to set service security)
id: 6ca3b6ff-dbdf-57ab-8bac-4d30f3a5af3b
status: experimental
description: 'Service registry key with ACL preventing standard enumeration tools from listing it - service exists but invisible to sc.exe and services.msc Query block: sc.exe used to set service security.'
references:
    - MITRE ATT&CK T1543.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1543.003
    - attack.g0096
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 1
        process.name: sc.exe
        process.command_line: '*sdset*'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1543.003
    technique_name: 'Create or Modify System Process: Windows Service'
    os: win
    query_block: 2 of 2
    original_query: |-
        winlog.event_id: 1
        AND process.name: "sc.exe"
        AND process.command_line: *sdset*
    sysmon: |-
        // EID 12 (registry key create/delete) or EID 13
        // catching changes to Service Security Descriptor:
        EventID=13
        TargetObject=*\Services\*\Security
        // The Security value contains the ACL for the
        // service registry key - modifications to it
        // indicate access control changes

        // Also watch sc.exe sdset commands:
        EventID=1
        Image=*\sc.exe
        CommandLine matches: *sdset*
        // sc sdset modifies a service's security descriptor
    powershell: |-
        # Hunt for hidden services - registry enumeration bypass
        # Method 1: Compare what sc.exe sees vs what's in the registry

        # What sc.exe / Get-Service can see:
        $visible = Get-Service | Select-Object -ExpandProperty Name

        # What's actually in the registry:
        $registry = Get-ChildItem "HKLM:\SYSTEM\CurrentControlSet\Services" |
          Where-Object { $_.GetValue('ImagePath') } |
          Select-Object -ExpandProperty PSChildName

        # Diff: services in registry but not visible to Get-Service
        $hidden = Compare-Object -ReferenceObject $registry -DifferenceObject $visible |
          Where-Object { $_.SideIndicator -eq '<=' } |
          Select-Object -ExpandProperty InputObject

        if ($hidden) {
          Write-Host "Potentially hidden services detected:" -ForegroundColor Red
          $hidden | ForEach-Object {
            $key = "HKLM:\SYSTEM\CurrentControlSet\Services\$_"
            $props = Get-ItemProperty $key -EA SilentlyContinue
            [PSCustomObject]@{
              Name = $_
              ImagePath = $props.ImagePath
              DisplayName = $props.DisplayName
              RegistryKey = $key
            }
          }
        } else {
          Write-Host "No hidden services detected via registry/SCM diff" -ForegroundColor Green
        }

        # Method 2: Check for modified service security descriptors
        # (requires reading the binary SD - shown conceptually)
        # Use: sc.exe sdshow <ServiceName>
        # Compare against known-good SD for that service type

        # Method 3: Enumerate service security from registry
        Get-ChildItem "HKLM:\SYSTEM\CurrentControlSet\Services" | ForEach-Object {
          $sec = (Get-ItemProperty $_.PSPath -Name Security -EA SilentlyContinue).Security
          if ($sec) {
            # Security value present means SD has been explicitly set
            # Default services don't have an explicit Security value -
            # they inherit from the parent key
            [PSCustomObject]@{
              ServiceName = $_.PSChildName
              HasExplicitSecurity = $true
              ImagePath = (Get-ItemProperty $_.PSPath -Name ImagePath -EA SilentlyContinue).ImagePath
            }
          }
        }
    registry: |-
        Service hiding via ACL manipulation:

        How it works:
        - Service registry key has a SECURITY_DESCRIPTOR (SD)
          controlling who can read/modify the key
        - Standard service SDs allow SYSTEM, Administrators,
          and authenticated users to enumerate
        - Adversary modifies SD to DENY access to
          Administrators and SYSTEM (or just deny READ to
          the specific accounts used by sc.exe / services.msc)
        - sc.exe and services.msc can't see the service
        - Service still runs because the SCM uses a different
          privilege level than user enumeration tools

        The Security value:
        HKLM\SYSTEM\CurrentControlSet\Services\<Name>\Security
        - Binary REG_BINARY value containing the SD
        - Optional - if absent, service uses default SD
        - Presence of this value on a non-system service
          is itself unusual

        Detection approaches:

        1. Diff registry vs SCM enumeration (shown in script):
           - List services from HKLM\SYSTEM\...Services\
           - List services from Get-Service / sc query
           - Difference = services hidden via ACL

        2. Inspect service SD directly:
           sc.exe sdshow <ServiceName>
           - Returns SDDL string
           - Compare against known-good SD format
           - Adversary SDs typically contain explicit DENY
             ACEs (D:...DENY... in SDDL format)

        3. Check for explicit Security value:
           - Most services don't have an explicit Security
             value (they inherit from parent)
           - A non-system service with explicit Security
             value warrants investigation

        Reset hidden service ACL (if found):
        sc.exe sdset <ServiceName>
          D:(A;;CCLCSWRPWPDTLOCRRC;;;SY)
          (A;;CCDCLCSWRPWPDTLOCRSDRCWDWO;;;BA)
          (A;;CCLCSWLOCRRC;;;IU)
          (A;;CCLCSWLOCRRC;;;SU)
        (restores default service SD)

        Memory forensics for confirmation:
        - Volatility psscan: lists running processes
          including those hidden from standard tools
        - Match running svchost / service processes
          against registry-enumerated services to identify
          processes belonging to hidden services
    ossdetect: |-
        Sigma:
        - registry_event_service_security_descriptor_change.yml
        - proc_creation_win_sc_sdset.yml

        Atomic Red Team:
        - T1543.003 has tests for service modification
        - Hidden service variant less commonly automated

        Hayabusa:
        - ServiceSDModification rules (EID 13)
        - ScSdSetExecution rules (EID 1)

        Velociraptor:
        - Windows.System.Services
          (enumerates from registry - surfaces hidden services
          that sc.exe / Get-Service cannot see)
        - Windows.Registry.NTUser
        - Best tool for fleet-wide hidden service hunting

        Reference:
        - Microsoft service security documentation
        - Windows Internals 7th edition - chapter on services
        - Various APT analyses describing service hiding
          tradecraft

        Manual workflow:
        1. Run registry/SCM diff script (shown above)
        2. For any flagged services, run sc.exe sdshow
        3. Compare SDDL against baseline
        4. Examine ImagePath and binary signature
        5. Restore default ACL only after collecting
           forensic evidence
    tools: |-
        Service hiding tradecraft:

        Known adversary toolkits with hidden service support:
        - Various custom China-nexus implants
        - Some advanced ransomware variants
        - Custom red team implants for evasion

        Standard implementation:
        - After installing the service, run:
          sc.exe sdset <ServiceName> <restrictive_SDDL>
        - Or modify the Security value in registry directly
        - Tools available for SDDL manipulation:
          - sc.exe sdshow / sdset (built-in)
          - PsService (Sysinternals)
          - Custom code using Win32 API
            (RegSetKeySecurity, SetServiceObjectSecurity)

        Why this technique is relatively rare:
        - More sophisticated than typical malware
        - Many adversaries don't bother because basic
          persistence works fine against most defenders
        - When you DO see it, treat as high-confidence
          evidence of capable / deliberate operator

        Detection note:
        - Defenders running sc query or services.msc will
          miss hidden services entirely
        - The registry diff approach is the most reliable
          detection - works even against SDs that deny
          access to all standard enumeration paths
        - Velociraptor's service enumeration uses registry
          reads, so it surfaces hidden services that sc.exe
          cannot see

        Real-world examples:
        - Reported in Mandiant and CrowdStrike intrusion
          reports on advanced operators
        - Less common than ServiceDll hijack but appears
          in long-dwell espionage cases
    actors:
        - APT41
        - APT41
    analyst_notes: 'Service hiding via security descriptor manipulation is one of the more sophisticated persistence techniques in Windows - it specifically targets the gap between how services are stored (registry) and how they''re typically enumerated (Service Control Manager via sc.exe). The technique works because Windows allows arbitrary ACLs on service registry keys, and the standard enumeration tools respect those ACLs - so if the ACL denies your account read access, the service appears not to exist from your perspective. The registry-vs-SCM diff approach is the most reliable detection because it bypasses the ACL entirely: you''re reading the registry directly and comparing against what the SCM returns. Any service in the registry that doesn''t appear in the SCM enumeration is either hidden, broken, or in an unusual state - all warrant investigation. This is a relatively rare technique - more common in APT-level operations than commodity malware - so finding it is a strong signal of deliberate compromise. Pair with the ServiceDll hijack indicator: an adversary using ServiceDll hijack AND hiding the service via SD manipulation has built a particularly stealthy persistence layer that defeats most standard defender workflows. Velociraptor''s Windows.System.Services artifact reads directly from the registry and is the best fleet-wide tool for surfacing this technique.'
