# sigma/rules/host/persistence/t1543_003_service_inventory_with_suspicious_configuration_2.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1543.003 - Service Inventory with Suspicious Configuration (Filter to suspicious service file paths)
id: 30cd4cb8-73c8-5870-ad48-c4cf0b9f120f
status: experimental
description: 'Enumerate all registered services and flag those with non-standard paths, unsigned binaries, blank descriptions, or other anomalous configuration Query block: Filter to suspicious service file paths.'
references:
    - MITRE ATT&CK T1543.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1543.003
    - attack.g0016
    - attack.g0096
    - attack.g0032
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 7045
        winlog.event_data.ImagePath:
            - '*\AppData\*'
            - '*\Temp\*'
            - '*\ProgramData\*'
            - '*\Public\*'
            - '*powershell*'
            - '*cmd.exe*'
            - '*rundll32*'
            - '*mshta*'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1543.003
    technique_name: 'Create or Modify System Process: Windows Service'
    os: win
    query_block: 3 of 3
    original_query: |-
        winlog.event_id: 7045
        AND winlog.event_data.ImagePath: (*\AppData\* OR *\Temp\* OR *\ProgramData\* OR *\Public\* OR *powershell* OR *cmd.exe* OR *rundll32* OR *mshta*)
    sysmon: |-
        // This indicator is hunting-oriented - PowerShell
        // registry sweep, not Sysmon detection.
        //
        // Sysmon detects service CREATION in real-time via
        // EID 13 (Run/Services registry write) and Windows
        // System Event 7045 - both covered in T1569.002.
        // This indicator finds services that were established
        // before logging coverage began.
        //
        // Sysmon contribution: re-check service event history:

        EventID=13 (registry value set)
        TargetObject=*\Services\*\ImagePath
        // Useful for retroactive timeline of when a
        // service was installed or modified

        // Windows System Event 7045 (service installed):
        EventID=7045
        LogName=System
        // Captures all service installations including
        // those by sc.exe, PowerShell New-Service,
        // PsExec PSEXESVC, and direct API calls
    powershell: |-
        # Comprehensive service inventory hunt - flag suspicious config

        $results = @()

        # Pull all services with their full registry config
        Get-ChildItem "HKLM:\SYSTEM\CurrentControlSet\Services" | ForEach-Object {
          $svc = $_
          $props = Get-ItemProperty $svc.PSPath -ErrorAction SilentlyContinue
          if (-not $props.ImagePath) { return }

          $name = $svc.PSChildName
          $ipRaw = $props.ImagePath
          # Strip leading quotes, extract just the binary path
          $ipPath = $ipRaw.TrimStart('"') -replace '^([^"]+).*','$1'
          $ipPath = [System.Environment]::ExpandEnvironmentVariables($ipPath)

          $suspicious = @()

          # Check 1: ImagePath in user-writable directory
          if ($ipPath -match '(AppData|Temp|ProgramData|\\Users\\Public|\\Downloads)') {
            $suspicious += "ImagePath in user-writable path"
          }

          # Check 2: ImagePath outside standard system/program dirs
          $standardPaths = @('C:\\Windows\\','C:\\Program Files\\','C:\\Program Files \(x86\)\\')
          $inStandard = $false
          foreach ($p in $standardPaths) {
            if ($ipPath -match "^$p") { $inStandard = $true; break }
          }
          if (-not $inStandard) {
            $suspicious += "ImagePath outside standard directories"
          }

          # Check 3: ImagePath is an interpreter (LOLBin service)
          if ($ipPath -match '(powershell|cmd\.exe|wscript|cscript|mshta|rundll32|regsvr32)\.exe$') {
            $suspicious += "Service wraps interpreter (LOLBin pattern)"
          }

          # Check 4: Binary unsigned (check Authenticode)
          if (Test-Path $ipPath) {
            $sig = Get-AuthenticodeSignature $ipPath -ErrorAction SilentlyContinue
            if ($sig.Status -ne 'Valid') {
              $suspicious += "Binary unsigned or invalid signature ($($sig.Status))"
            }
          } else {
            $suspicious += "Binary file does not exist at ImagePath"
          }

          # Check 5: Blank Description and DisplayName (anomalous)
          if ([string]::IsNullOrWhiteSpace($props.Description) -and
              [string]::IsNullOrWhiteSpace($props.DisplayName)) {
            $suspicious += "Blank Description AND DisplayName"
          }

          # Check 6: Auto-start service with no DependOnService
          # (anomalous for non-Microsoft services)
          if ($props.Start -eq 2 -and -not $props.DependOnService -and
              $ipPath -notmatch 'C:\\Windows\\') {
            $suspicious += "Auto-start with no dependencies, non-system path"
          }

          if ($suspicious.Count -gt 0) {
            $results += [PSCustomObject]@{
              ServiceName = $name
              DisplayName = $props.DisplayName
              ImagePath = $ipRaw
              StartType = switch ($props.Start) {
                0 {'Boot'} 1 {'System'} 2 {'Auto'}
                3 {'Demand'} 4 {'Disabled'} default {$_}
              }
              RunAs = $props.ObjectName
              Description = $props.Description
              Suspicious = ($suspicious -join '; ')
            }
          }
        }

        $results | Sort-Object ServiceName | Format-Table -AutoSize
    registry: |-
        Service registry structure:
        HKLM\SYSTEM\CurrentControlSet\Services\<ServiceName>\
          ImagePath     - binary path + arguments
          Start         - 0=Boot, 1=System, 2=Auto, 3=Demand, 4=Disabled
          Type          - 0x10=Win32OwnProcess, 0x20=Win32ShareProcess (svchost),
                          0x110=interactive Win32OwnProcess
          ObjectName    - run-as account (LocalSystem, NetworkService, etc.)
          Description   - human-readable description
          DisplayName   - shown in services.msc
          DependOnService - service dependencies (REG_MULTI_SZ)
          ErrorControl  - error handling on failure

        Service hunting workflow:
        1. Run the PowerShell sweep above
        2. For each flagged service, investigate:
           - File hash of ImagePath binary (VirusTotal lookup)
           - When was the service key last modified?
             (Get-Item <key> | Select LastWriteTime)
           - Process activity history for that service name
           - Network connections from service process

        Baseline approach:
        - Capture service inventory on known-good reference
          system: Get-Service | Export-CSV baseline.csv
        - Diff suspect host:
          Compare-Object (Import-CSV baseline.csv)
            (Get-Service | Select Name,DisplayName)
        - New service entries warrant investigation

        Last-write timestamps for IR timelining:
        (Get-Item "HKLM:\SYSTEM\CurrentControlSet\Services\<name>").LastWriteTime
        - Tells you when the service was installed or last
          modified - useful for correlating with other
          intrusion timeline events

        Standard service binary locations:
        - C:\Windows\System32\
        - C:\Windows\SysWOW64\
        - C:\Program Files\<vendor>\
        - C:\Program Files (x86)\<vendor>\
        - Anything else = investigate
    ossdetect: |-
        Sigma:
        - Real-time detection rules cover creation
          (proc_creation_win_sc_service_creation_*.yml)
        - Hunting existing services is a manual workflow

        Atomic Red Team:
        - T1543.003 Test #1 (sc.exe service creation)
        - T1543.003 Test #2 (PowerShell New-Service)
        - Use to validate the hunt: install service via
          Atomic, run the PowerShell sweep, confirm it
          surfaces in your detection

        Hayabusa:
        - ServiceCreationSuspBinPath rules (real-time)
        - For hunting, use PowerShell or autoruns

        Velociraptor:
        - Windows.System.Services (best fleet-wide artifact)
        - Returns full service config including signature
          status for each ImagePath binary

        Reference - SANS DFIR poster:
        'Windows Services Forensics'
        - Comprehensive coverage of service abuse patterns
        - Annual updates with new techniques

        Microsoft documentation:
        - Service security and access rights
        - Useful for understanding the technique surface
    tools: |-
        Hunting / enumeration tools:

        PowerShell (the script above):
        - No tool deploy required
        - Works on any Windows
        - Best for ad-hoc IR triage

        Sysinternals autoruns.exe:
        - Services tab - all registered services
        - Highlights unsigned and unverified binaries
        - VirusTotal integration for reputation
        - Single best tool for manual triage

        Velociraptor:
        - Windows.System.Services
          (full service enumeration artifact)
        - Returns ImagePath, signature status, last-write
          time, and dependencies for each service
        - Best for fleet-wide service hunting

        KAPE (Eric Zimmerman):
        - Targets/Compound/SystemHive.tkape
          (collects SYSTEM registry hive for offline analysis)
        - Combined with RECmd plugins for service enumeration

        Sigma rules apply mostly to creation-time
        detection (covered in T1569.002) - hunting
        established services is more of a manual workflow

        Adversary tool patterns to recognize:
        - Cobalt Strike services: random short names,
          ImagePath often points to a DLL via rundll32,
          or to a small EXE in user-writable location
        - PSEXESVC: well-known PsExec pattern (C:\Windows\PSEXESVC.exe)
        - SharPersist: configurable but typically uses
          legitimate-sounding names with binaries in
          C:\ProgramData
        - Empire / Metasploit: predictable naming patterns
          (mostly addressed by their respective module docs)
    actors:
        - APT29
        - APT41
        - Lazarus
    analyst_notes: 'Service hunting is one of the most important persistence-detection workflows in Windows IR. The PowerShell sweep covers the high-fidelity flags: services with binaries in user-writable paths, services wrapping interpreters (LOLBin pattern), unsigned binaries, missing files (orphaned service entries), blank metadata. Each flag alone is moderate confidence; multiple flags on the same service is high confidence. The baseline-and-diff approach is especially powerful for services because legitimate service inventories are stable over time - the same set of services with the same configs runs on every workstation of the same build. Deviations from baseline are immediately suspicious. The unsigned binary check is particularly useful: legitimate Microsoft and major-vendor services are virtually always Authenticode-signed. An unsigned binary running as a service is rare in enterprise environments outside of in-house developed software. Pair this hunt with autoruns.exe for the highest coverage - autoruns includes services in its 200+ autostart inventory and adds VirusTotal reputation lookups that the PowerShell sweep doesn''t have. The last-write timestamp trick on the service registry key is genuinely useful in IR timelining: it tells you when the service was installed or last modified, which often correlates with initial compromise events visible in other logs.'
