# sigma/rules/host/persistence/t1543_003_servicedll_hijack_svchost_hosted_service_persist_1.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1543.003 - ServiceDLL Hijack (svchost-hosted Service Persistence) (svchost loading unsigned DLL from non-system path)
id: 211c12cb-1d80-5be7-92b8-ecb54a3d774d
status: experimental
description: 'Service of Type=0x20 (Win32ShareProcess) with ServiceDll registry value pointing to user-writable or unsigned DLL - svchost-loaded malicious DLL Query block: svchost loading unsigned DLL from non-system path.'
references:
    - MITRE ATT&CK T1543.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1543.003
    - attack.g0096
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 7
        process.name: svchost.exe
        file.code_signature.signed: 'false'
    filter_0:
        file.path:
            - C:\\Windows\\System32\\*
            - C:\\Windows\\SysWOW64\\*
    condition: selection and not filter_0
falsepositives:
    - legitimate Microsoft services use ServiceDll extensively (it's the standard svchost service pattern), so the filter must include 'path outside System32/SysWOW64' and 'unsigned' to surface anomalies
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1543.003
    technique_name: 'Create or Modify System Process: Windows Service'
    os: win
    query_block: 2 of 2
    original_query: |-
        winlog.event_id: 7
        AND process.name: "svchost.exe"
        AND file.code_signature.signed: false
        AND NOT file.path: ("C:\\Windows\\System32\\*" OR "C:\\Windows\\SysWOW64\\*")
    sysmon: |-
        // Real-time detection - registry write to ServiceDll:
        EventID=13
        TargetObject=*\Services\*\Parameters\ServiceDll
        Details matches:
          *\AppData\* OR *\Temp\*
          OR *\ProgramData\* OR *\Users\Public\*
          OR <DLL path outside System32>

        // Real-time detection - svchost loading unsigned DLL:
        EventID=7 (Image Load)
        Image=*\svchost.exe
        ImageLoaded NOT in standard system paths
        Signed=false
    powershell: |-
        # Hunt for ServiceDll values pointing to suspicious paths
        Get-ChildItem "HKLM:\SYSTEM\CurrentControlSet\Services" | ForEach-Object {
          $svc = $_
          $paramKey = Join-Path $svc.PSPath 'Parameters'
          if (Test-Path $paramKey) {
            $dll = (Get-ItemProperty $paramKey -Name ServiceDll -EA SilentlyContinue).ServiceDll
            if ($dll) {
              $dllExpanded = [System.Environment]::ExpandEnvironmentVariables($dll)

              $suspicious = @()

              # User-writable path
              if ($dllExpanded -match '(AppData|Temp|ProgramData|\\Users\\Public)') {
                $suspicious += "ServiceDll in user-writable path"
              }

              # Outside standard system / program directories
              if ($dllExpanded -notmatch '^(C:\\Windows\\|C:\\Program Files)') {
                $suspicious += "ServiceDll outside standard paths"
              }

              # Signature check
              if (Test-Path $dllExpanded) {
                $sig = Get-AuthenticodeSignature $dllExpanded -EA SilentlyContinue
                if ($sig.Status -ne 'Valid') {
                  $suspicious += "ServiceDll unsigned ($($sig.Status))"
                }
              } else {
                $suspicious += "ServiceDll file does not exist"
              }

              if ($suspicious.Count -gt 0) {
                [PSCustomObject]@{
                  ServiceName = $svc.PSChildName
                  ServiceDll = $dll
                  Resolved = $dllExpanded
                  Suspicious = ($suspicious -join '; ')
                }
              }
            }
          }
        }

        # Cross-reference: which svchost group hosts this service?
        # Get-ItemProperty "HKLM:\SYSTEM\CurrentControlSet\Services\<name>" -Name ImagePath
        # The ImagePath typically contains: %SystemRoot%\system32\svchost.exe -k <group>
    registry: |-
        svchost service architecture:
        - svchost.exe hosts multiple services in one process
        - Service "groups" allow related services to share
          a single svchost instance (memory efficiency)
        - Each service has:
          HKLM\SYSTEM\CurrentControlSet\Services\<Name>\
            ImagePath = %SystemRoot%\system32\svchost.exe
                        -k <group_name>
            Parameters\
              ServiceDll = <path to DLL implementing service>
              ServiceMain = <export function name, optional>

        The hijack technique:
        - Adversary creates a service of Type=0x20 (shared)
          with ImagePath pointing to legitimate svchost.exe
        - Sets Parameters\ServiceDll to point to attacker-
          controlled DLL
        - When service starts, svchost.exe (trusted Microsoft
          binary) loads the attacker DLL
        - Process tree shows svchost.exe as the running
          process - looks legitimate to most monitoring

        Why adversaries favor this:
        - Hides under svchost.exe process name (40-100 svchosts
          run on a typical Windows host - blending is easy)
        - ImagePath is the legitimate Microsoft svchost.exe -
          bypasses path-based detection
        - Detection requires looking at the ServiceDll parameter
          rather than the ImagePath
        - The malicious DLL runs in a SYSTEM-privileged
          process context

        Hunting workflow:
        - Enumerate Services with ImagePath=svchost.exe -k *
        - For each, check Parameters\ServiceDll
        - Flag ServiceDll values outside System32/SysWOW64

        Defensive registry:
        - KnownDLLs (HKLM\SYSTEM\CurrentControlSet\Control\
          Session Manager\KnownDLLs) does NOT protect against
          ServiceDll hijack since the DLL is loaded by explicit
          path, not by search order
    ossdetect: |-
        Sigma:
        - registry_event_servicedll_susp_path.yml
        - image_load_win_svchost_susp_dll.yml
        - registry_event_servicedll_modification.yml

        Atomic Red Team:
        - T1543.003 (service variants)
        - Some tests cover ServiceDll modification

        Hayabusa:
        - ServiceDllSuspPath rules (EID 13)
        - SvchostUnsignedDllLoad detection (EID 7)

        Velociraptor:
        - Windows.System.Services (returns Parameters\ServiceDll
          alongside ImagePath for each service)
        - Windows.System.DLL (loaded DLLs per process)

        Reference - Microsoft documentation:
        - 'Working with the AllJoyn Service Manager'
          and similar svchost service implementation docs
        - Useful for understanding the legitimate side
          of the technique surface

        Threat hunting references:
        - PlugX malware family analyses
        - Various China-nexus APT reports documenting
          ServiceDll hijack tradecraft
    tools: |-
        Adversary frameworks using ServiceDll hijack:
        - Various long-dwell espionage operators
        - ShadowPad / PlugX variants
        - Some commodity malware (less common - prefers
          simpler Run key persistence)

        Build process:
        - Adversary writes a DLL exporting ServiceMain function
        - Optionally exports SvchostPushServiceGlobals
          (some service groups require it)
        - DLL placed in attacker-chosen location
        - Service registered via sc.exe or registry edit

        Reference DLL structure:
        - DllMain entry point
        - ServiceMain export (the service implementation)
        - ServiceCtrlHandler for SCM control messages
        - Most malicious ServiceDlls implement minimal
          legitimate-looking service interface, then run
          attacker code in a separate thread

        Detection note: this is a more sophisticated
        persistence than simple service installation,
        typically associated with deliberate operators
        rather than commodity malware. Finding ServiceDll
        hijack on a host is usually a strong indicator
        of APT-level compromise rather than opportunistic
        infection.

        Service group context:
        - netsvcs (the most populated group - 30+ services)
        - LocalService, NetworkService, LocalSystemNetworkRestricted
        - Each group has its own svchost instance
        - Adversaries often target netsvcs because it's
          the most crowded and most generic
    actors:
        - APT41
        - ShadowPad
        - PlugX
        - APT41
    analyst_notes: 'ServiceDll hijack is a more sophisticated persistence technique than basic service installation - it specifically targets the svchost-hosted service model where multiple services share a single host process. The detection pivot is that you''re not looking at the ImagePath (which is legitimate svchost.exe), you''re looking at the Parameters\ServiceDll value which actually points to the executable code. This is a high-value indicator when found because it''s predominantly associated with deliberate, capable operators rather than commodity malware - most opportunistic infections use simpler Run key or service-with-EXE persistence rather than going through the trouble of building a proper ServiceDll. ShadowPad and PlugX malware families are the canonical examples and are heavily associated with China-nexus APT operations. The KnownDLLs registry protection does NOT apply here since the DLL is loaded by explicit path, not search order - so this technique works regardless of DLL search order hardening. False positives: legitimate Microsoft services use ServiceDll extensively (it''s the standard svchost service pattern), so the filter must include ''path outside System32/SysWOW64'' and ''unsigned'' to surface anomalies. Pair the Sysmon EID 13 real-time detection with the PowerShell hunt for full coverage.'
