# sigma/rules/host/persistence/t1547_001_real_time_registry_run_key_modification_tighten.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1547.001 - Real-Time Registry Run Key Modification (Tighten by suspicious target path / interpreter)
id: 64387d73-805b-5abe-a5a0-4b9d9125dbc6
status: experimental
description: 'Sysmon EID 13 - registry value written to Run/RunOnce key, especially with target in user-writable path Query block: Tighten by suspicious target path / interpreter.'
references:
    - MITRE ATT&CK T1547.001
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1547.001
    - attack.g0016
    - attack.g0096
    - attack.g0059
    - attack.g0032
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        winlog.event_id: 13
        registry.path:
            - '*\Run\*'
            - '*\RunOnce\*'
        registry.data.strings:
            - '*\AppData\*'
            - '*\Temp\*'
            - '*\ProgramData\*'
            - '*\Public\*'
            - '*powershell*'
            - '*cmd.exe*'
            - '*wscript*'
            - '*mshta*'
            - '*rundll32*'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1547.001
    technique_name: 'Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder'
    os: win
    query_block: 2 of 3
    original_query: |-
        winlog.event_id: 13
        AND registry.path: (*\Run\* OR *\RunOnce\*)
        AND registry.data.strings: (*\AppData\* OR *\Temp\* OR *\ProgramData\* OR *\Public\* OR *powershell* OR *cmd.exe* OR *wscript* OR *mshta* OR *rundll32*)
    sysmon: |-
        EventID=13
        TargetObject matches (any of):
          *\Microsoft\Windows\CurrentVersion\Run\*
          *\Microsoft\Windows\CurrentVersion\RunOnce\*
          *\Microsoft\Windows\CurrentVersion\RunServices\*
          *\Microsoft\Windows\CurrentVersion\RunServicesOnce\*
          *\Microsoft\Windows\CurrentVersion\Explorer\
            User Shell Folders\Startup
          *\Microsoft\Windows NT\CurrentVersion\Winlogon\
            Userinit
          *\Microsoft\Windows NT\CurrentVersion\Winlogon\Shell
          *\Microsoft\Windows NT\CurrentVersion\Image File
            Execution Options\*\Debugger

        // Value being written (the Details field) usually
        // contains the binary path - check for user-writable
        // paths or interpreter invocations:
        Details matches:
          *\AppData\* OR *\Temp\* OR *\ProgramData\*
          OR *\Users\Public\*
          OR *powershell* OR *cmd.exe* OR *wscript*
          OR *mshta* OR *rundll32* OR *regsvr32*
    powershell: |-
        # Hunt for real-time Run key modifications (Sysmon EID 13)
        $autoStartKeys = @(
          '\Run\','\RunOnce\','\RunServices\','\RunServicesOnce\',
          'Userinit','Winlogon\Shell',
          'Image File Execution Options\.*\Debugger'
        )

        Get-WinEvent -FilterHashtable @{
          LogName='Microsoft-Windows-Sysmon/Operational';
          ID=13
        } | Where-Object {
          $target = $_.Properties[4].Value
          $autoStartKeys | Where-Object { $target -match $_ }
        } | Select TimeCreated,
          @{n='TargetKey';e={$_.Properties[4].Value}},
          @{n='Details';e={$_.Properties[5].Value}},
          @{n='Image';e={($_.Properties[3].Value -split '\\')[-1]}},
          @{n='User';e={$_.Properties[2].Value}} |
          Sort-Object TimeCreated -Descending

        # Startup folder file creates (Sysmon EID 11)
        Get-WinEvent -FilterHashtable @{
          LogName='Microsoft-Windows-Sysmon/Operational';
          ID=11
        } | Where-Object {
          $_.Properties[0].Value -match 'Start Menu\\Programs\\Start[Uu]p'
        } | Select TimeCreated,
          @{n='File';e={$_.Properties[0].Value}},
          @{n='CreatingProcess';e={$_.Properties[5].Value}}
    registry: |-
        Primary autostart registry locations:

        User-context (no admin required):
        HKCU\Software\Microsoft\Windows\CurrentVersion\Run
        HKCU\Software\Microsoft\Windows\CurrentVersion\RunOnce
        HKCU\Software\Microsoft\Windows NT\CurrentVersion\
          Windows\Load
        HKCU\Software\Microsoft\Windows NT\CurrentVersion\
          Windows\Run

        Machine-context (requires admin):
        HKLM\SOFTWARE\Microsoft\Windows\CurrentVersion\Run
        HKLM\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce
        HKLM\SOFTWARE\Microsoft\Windows\CurrentVersion\
          RunServices
        HKLM\SOFTWARE\Microsoft\Windows\CurrentVersion\
          RunServicesOnce
        HKLM\SOFTWARE\Wow6432Node\Microsoft\Windows\
          CurrentVersion\Run (32-bit equivalent)

        Winlogon variants:
        HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\
          Winlogon\Userinit
        HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\
          Winlogon\Shell

        Image File Execution Options (IFEO debugger hijack):
        HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\
          Image File Execution Options\<TargetEXE>\Debugger
          - Value here causes Debugger to launch instead of
            the target EXE when that EXE is invoked
          - Classic abuse: hijack a frequently-run binary

        Startup folder file locations:
        - All-users startup:
          C:\ProgramData\Microsoft\Windows\Start Menu\
            Programs\StartUp\
        - Per-user startup:
          %APPDATA%\Microsoft\Windows\Start Menu\
            Programs\Startup\
        - Anything (.exe/.lnk/.bat/.vbs) here runs at logon
    ossdetect: |-
        Sigma:
        - registry_event_run_key_modification.yml
        - registry_event_run_key_susp_path.yml
        - registry_event_winlogon_userinit_susp_change.yml
        - registry_event_ifeo_debugger_set.yml
        - file_event_win_startup_folder_persistence.yml

        Atomic Red Team:
        - T1547.001 (multiple Run key variants tested)
        - T1547.001 Test #1 (Run key add via reg.exe)
        - T1547.001 Test #5 (Startup folder file drop)

        Hayabusa:
        - RegistryRunKeyMod rules (EID 13)
        - StartupFolderFileWrite rules (EID 11)
        - High coverage of the autostart locations

        Velociraptor:
        - Windows.Sys.AutoRuns (the dedicated artifact)
        - Windows.Forensics.Autoruns
        - Windows.Registry.NTUser (HKCU enumeration)
        - Best fleet-wide tooling for this technique

        Sysinternals autoruns.exe:
        - The canonical persistence-hunting tool
        - Enumerates 200+ autostart locations
        - Highlights unsigned entries in red/yellow
        - VirusTotal integration for reputation checking
        - For manual IR triage - close to essential
    tools: |-
        Almost universal across malware families:
        - Commodity malware (Emotet, TrickBot, IcedID,
          QakBot, Agent Tesla, Formbook, etc.)
        - Most RATs (NjRAT, DarkComet, Remcos, etc.)
        - Most ransomware staging
        - Manual operator persistence
        - Some legitimate vendor software (false positive source)

        Tool / framework support:
        - SharPersist (Brian Mulligan):
          Built-in support for all Run key variants
          Standard red team persistence tool
        - Empire / Starkiller persistence modules
        - Cobalt Strike persistence aggressor scripts
        - Metasploit persistence_exe / registry_persistence
        - Built-in Windows reg.exe and PowerShell:
          reg add "HKCU\Software\Microsoft\Windows\
            CurrentVersion\Run" /v Updater /t REG_SZ
            /d "C:\Users\victim\AppData\evil.exe"

        Common adversary value name patterns:
        - Spoofed Microsoft names: "Windows Update",
          "Microsoft Defender", "Adobe Updater"
        - Random strings: "svc_a8f3", "winupd_x64"
        - Blank or whitespace value names
        - Names mimicking installed legitimate software
    actors:
        - APT29
        - APT41
        - APT35
        - Lazarus
    analyst_notes: 'Registry Run Keys + Startup Folders is the single most common persistence technique on Windows - if a piece of malware persists, there''s a good chance it''s using one of these mechanisms. The technique is so common that detection cannot be ''alert on every Run key modification'' - that would be overwhelming. Tune by value-data context: alert on Run key writes where the target path is in AppData/Temp/ProgramData, contains an interpreter (powershell/cmd/wscript/mshta), or matches known suspicious patterns. HKCU variants don''t require admin privileges - any compromised user can persist their own session through HKCU\Run. The IFEO Debugger trick is worth understanding: setting an Image File Execution Options Debugger value causes Windows to launch the Debugger value instead of the original binary when the original is invoked - so hijacking sethc.exe (Sticky Keys), magnify.exe, or osk.exe creates accessibility-feature persistence that triggers from the logon screen. Pair this real-time detection with the next indicator (hunting all locations) for full coverage - real-time catches new persistence as it''s installed; hunting catches persistence that was installed before logging started or that bypassed Sysmon coverage.'
