# sigma/rules/host/persistence/t1543_002_malicious_service_unit_file_drop_and_enable_new.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1543.002 - Malicious .service Unit File Drop and Enable (New .service or .timer file creation in systemd directories)
id: 2c3bd4e6-6e50-5f44-844c-216c35125080
status: experimental
description: 'New .service file written outside package manager into /etc/systemd/system/, /usr/lib/systemd/system/, or ~/.config/systemd/user/, followed by systemctl enable or daemon-reload Query block: New .service or .timer file creation in systemd directories. NOTE: partial conversion - one or more source clauses could not be expressed in stock Sigma. See tonk.dropped_clauses.'
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1543.002
    - attack.g0096
    - attack.g0010
    - attack.g0016
    - attack.g0094
logsource:
    category: file_event
    product: windows
detection:
    selection:
        event.module: file_integrity
        file.path:
            - /etc/systemd/system/*
            - /usr/lib/systemd/system/*
            - /usr/local/lib/systemd/system/*
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: partial
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1543.002
    technique_name: 'Create or Modify System Process: Systemd Service'
    os: linux
    query_block: 1 of 5
    original_query: |-
        event.module: "file_integrity"
        AND file.path: (/etc/systemd/system/* OR /usr/lib/systemd/system/* OR /usr/local/lib/systemd/system/*)
        AND (file.path: *.service OR file.path: *.timer)
    sysmon: |-
        # Auditd rules for systemd service file persistence

        # System-level unit file directories (root required)
        -w /etc/systemd/system -p wa -k systemd_unit_write
        -w /usr/lib/systemd/system -p wa -k systemd_unit_write
        -w /usr/local/lib/systemd/system -p wa -k systemd_unit_write

        # Drop-in directories (common adversary bypass technique)
        -w /etc/systemd/system.conf.d -p wa -k systemd_unit_write

        # systemctl binary invocations (enable/disable/start)
        -a always,exit -F arch=b64 -S execve 
          -F path=/bin/systemctl -k systemctl_exec
        -a always,exit -F arch=b64 -S execve 
          -F path=/usr/bin/systemctl -k systemctl_exec

        # systemd generator directories (advanced persistence)
        -w /etc/systemd/system-generators -p wa -k systemd_generator
        -w /usr/lib/systemd/system-generators -p wa -k systemd_generator
        -w /usr/local/lib/systemd/system-generators -p wa -k systemd_generator

        # Sysmon for Linux EID 11 (file create) — catches .service drops
        # Filter: TargetFilename contains ".service" or ".timer"
    powershell: |-
        #!/bin/bash
        # T1543.002 - Systemd Service Persistence Hunt

        echo "[*] === All active/failed services (non-package) ==="
        systemctl list-units --type=service --all --no-pager 2>/dev/null | 
          grep -v "^UNIT"

        echo ""
        echo "[*] === Services enabled for boot ==="
        systemctl list-unit-files --type=service --state=enabled --no-pager 2>/dev/null

        echo ""
        echo "[*] === Unit files NOT from package manager (manual drops) ==="
        # Installed pkg units are in /usr/lib/systemd; manual ones in /etc/systemd/system
        # Flag /etc/systemd/system/*.service that aren't symlinks to /usr/lib
        for f in /etc/systemd/system/*.service; do
          [ -f "$f" ] || continue
          if ! readlink -f "$f" | grep -q "^/usr/lib/systemd"; then
            echo "[FLAG] Non-packaged service: $f"
            cat "$f"
            echo ""
          fi
        done

        echo ""
        echo "[*] === ExecStart paths in suspicious locations ==="
        grep -r "ExecStart" /etc/systemd/system/ /usr/lib/systemd/system/ 2>/dev/null | 
          grep -E "(/tmp/|/dev/shm/|/var/tmp/|/home/[^/]+/.[^/]+|/usr/local/sbin/.)" 

        echo ""
        echo "[*] === User-level systemd units (unprivileged persistence) ==="
        find /home -path "*/.config/systemd/user/*.service" -type f 2>/dev/null | 
          while read f; do
            echo "--- $f ---"
            cat "$f"
          done
        find /root/.config/systemd/user/ -name "*.service" 2>/dev/null | 
          while read f; do
            echo "--- $f ---"
            cat "$f"
          done

        echo ""
        echo "[*] === Recently modified unit files (7 days) ==="
        find /etc/systemd /usr/lib/systemd -name "*.service" -mtime -7 -ls 2>/dev/null

        echo ""
        echo "[*] === Systemd generators (advanced persistence) ==="
        ls -la /etc/systemd/system-generators/ 
               /usr/lib/systemd/system-generators/ 
               /usr/local/lib/systemd/system-generators/ 2>/dev/null
    registry: |-
        Systemd unit file persistence locations:

        System-level (root required):
          /etc/systemd/system/<name>.service     - admin-installed units
          /usr/lib/systemd/system/<name>.service - package-installed units
          /usr/local/lib/systemd/system/         - locally compiled packages
          /lib/systemd/system/                   - symlink to /usr/lib (Debian)

        Drop-in override directories:
          /etc/systemd/system/<name>.service.d/  - override fragments
          (can override ExecStart without touching the original unit)

        User-level (no root needed):
          ~/.config/systemd/user/<name>.service  - per-user unit
          /usr/lib/systemd/user/<name>.service   - user unit from package

        Systemd generator directories (advanced):
          /etc/systemd/system-generators/
          /usr/lib/systemd/system-generators/
          /usr/local/lib/systemd/system-generators/
          (executables here run at boot to dynamically create units)

        Enable persistence artifacts:
          /etc/systemd/system/multi-user.target.wants/<name>.service
          (symlink created by systemctl enable)

        Critical unit file directives to check:
          ExecStart=   - the binary or command being run
          User=        - if set to root on non-system service (flag)
          Restart=always  - ensures restart on death (persistence signal)
          [Install] WantedBy=multi-user.target  - boot enable
    ossdetect: |-
        Sigma rules:
        - file_event_lnx_systemd_service_file_creation.yml
        - proc_creation_lnx_systemctl_susp_service_enable.yml
        - file_event_lnx_systemd_generator_creation.yml
        - proc_creation_lnx_systemd_service_susp_path.yml

        Elastic detection rules:
        - Potential Persistence via Systemd Service Creation
        - Systemd Unit File Created
        - Systemd Service Started/Enabled via Suspicious Argument

        Atomic Red Team:
        - T1543.002 Test #1 (service unit file creation)
        - T1543.002 Test #2 (systemctl enable)
        - T1543.002 Test #3 (user-level unit file)

        PANIX:
        - --systemd-service flag, tests both root and user units
          github.com/Aegrah/PANIX

        Velociraptor:
        - Linux.Sys.Services (enumerate all systemd units)
        - Linux.Detection.Systemd (anomaly hunt)

        osquery:
        - SELECT * FROM systemd_units WHERE active_state = 'active';
        - SELECT * FROM systemd_units WHERE unit_file_state = 'enabled'
          AND source_path LIKE '/etc/systemd/system/%';

        Pepe Berba blog series:
        - 'Hunting for Persistence in Linux (Part 5): Systemd Generators'
        - Required reading for generator-based persistence
    tools: |-
        APT actors and malware using systemd persistence:

        Rocke (CN) - cryptominer; installs systemd services named after
          legitimate daemons (sshd-sync, system-update, irqbalanced)
          to run mining payloads with restart=always

        TeamTNT (mul) - cloud/container threat actor; drops .service
          files that start C2 agents and miners on compromised hosts

        UNC3524 (likely CN) - used systemd service persistence on
          compromised Linux servers during long-dwell espionage;
          Mandiant documented 18+ month dwell time

        Gomir (KP/Kimsuky) - Linux backdoor dropped as systemd service
          with ExecStart pointing to /usr/bin/timedatectl (spoofed path)

        Ebury (RU-nexus) - SSH backdoor rootkit also establishes
          systemd service for persistence on OpenSSH servers

        PANIX (simulation tool):
          ./panix.sh --systemd-service  - drops test .service file
          github.com/Aegrah/PANIX

        Common red team/attacker patterns:
          [Unit]
          Description=System Logger Service
          [Service]
          ExecStart=/usr/local/bin/.svc-daemon
          Restart=always
          [Install]
          WantedBy=multi-user.target

          systemctl enable --now malicious.service
          systemctl daemon-reload
    dropped_clauses:
        - 'unparsed: (file.path: *.service OR file.path: *.timer)'
    actors:
        - APT41
        - Turla
        - Rocke
        - UNC3524
        - Kimsuky
        - Ebury
        - TeamTNT
    analyst_notes: 'Systemd service persistence is the Linux equivalent of creating a Windows service - it requires root for system-level units but only a valid user account for user-level units in ~/.config/systemd/user. The asymmetry matters: system services run as root and start at boot; user services run as the owning user and start at user login. Attackers frequently use system services for the stronger persistence guarantee. The ExecStart path is the primary IOC: legitimate system services almost never point to /tmp, /dev/shm, or home directories. The Restart=always directive is a second-order IOC - it ensures the payload respawns if killed, exactly what a defender would do during remediation. For detection, the strongest signal is a new .service file appearing in /etc/systemd/system that was not created by a package manager (apt, yum, dpkg, rpm). Diff the mtime of all .service files against the last known package installation timestamp. Systemd generators are an advanced and underdetected variant: executables placed in generator directories run at boot before normal services and can dynamically create additional unit files, making them a persistent foothold that survives even removal of the visible .service artifact.'
