# sigma/rules/host/persistence/t1547_001_suspicious_value_names_and_spoofing_patterns_pur.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1547.001 - Suspicious Value Names and Spoofing Patterns (Pure-random-looking value names (regex))
id: a9ddeca7-b7b7-5c9a-abb2-36cabaf6b31d
status: experimental
description: 'Autostart entry value name mimics legitimate software, contains random strings, or has anomalous characteristics suggesting impersonation Query block: Pure-random-looking value names (regex).'
references:
    - MITRE ATT&CK T1547.001
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1547.001
    - attack.g0032
    - attack.g0096
    - attack.g0059
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        winlog.event_id: 13
        registry.path: '*\Run\*'
        registry.value|re: ^[a-z0-9]{8,16}$
    condition: selection
falsepositives:
    - False positives are inherent here - run this as a triage prioritizer, not a final alert
level: high
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1547.001
    technique_name: 'Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder'
    os: win
    query_block: 2 of 2
    original_query: |-
        winlog.event_id: 13
        AND registry.path: *\Run\*
        AND registry.value: /^[a-z0-9]{8,16}$/
    sysmon: |-
        // EID 13 with value-name-specific filtering:
        EventID=13
        TargetObject matches: *\Run\*

        // Suspicious value name patterns:
        // - Spoofing Microsoft/Adobe/Google product names:
        ValueName matches:
          *Microsoft* OR *Windows* OR *Defender*
          OR *Update* OR *Adobe* OR *Google*
          OR *Chrome* OR *Firefox*
        // AND target path is NOT in expected location:
        Details NOT matching:
          C:\Program Files\<expected_publisher>\*
          C:\Program Files (x86)\<expected_publisher>\*

        // - Random alphanumeric (machine-generated names):
        ValueName matches: [a-z0-9]{8,16} (regex - 8+ random chars)

        // - Blank or whitespace-only value names:
        ValueName="" OR ValueName=" " OR ValueName matches: ^\s+$
    powershell: |-
        # Hunt for spoofed / suspicious value names in Run keys
        $runKeys = @(
          'HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\Run',
          'HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce',
          'HKCU:\SOFTWARE\Microsoft\Windows\CurrentVersion\Run',
          'HKCU:\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce'
        )

        # Spoofable product names (commonly impersonated)
        $spoofableNames = 'Microsoft|Windows|Defender|Update|Adobe|Google|Chrome|Firefox|Apple|iTunes|Skype|Office|OneDrive|Dropbox|Java'

        # Expected publishers for each spoofed name
        $expectedPaths = @{
          'Microsoft' = 'Program Files\\(Microsoft|Common Files)'
          'Windows'   = 'C:\\Windows\\'
          'Adobe'     = 'Program Files\\Adobe'
          'Google'    = 'Program Files\\Google'
          'Chrome'    = 'Program Files\\Google\\Chrome'
        }

        foreach ($key in $runKeys) {
          if (Test-Path $key) {
            $props = Get-ItemProperty $key
            $props.PSObject.Properties | Where-Object {
              $_.Name -notmatch '^PS' -and $_.Value
            } | ForEach-Object {
              $name = $_.Name
              $value = $_.Value -as [string]
              $suspicious = @()

              # Check 1: spoofed name with wrong path
              if ($name -match $spoofableNames) {
                $matchedTerm = $matches[0]
                $expected = $expectedPaths[$matchedTerm]
                if ($expected -and $value -notmatch $expected) {
                  $suspicious += "Spoof: name=$matchedTerm but path doesn't match expected $expected"
                }
              }

              # Check 2: random-looking value name
              if ($name -match '^[a-z0-9]{8,16}$') {
                $suspicious += "Random value name: $name"
              }

              # Check 3: blank/whitespace name
              if ([string]::IsNullOrWhiteSpace($name) -or $name -match '^\s+$') {
                $suspicious += "Blank or whitespace value name"
              }

              # Check 4: value path in user-writable directory
              if ($value -match '(AppData|Temp|ProgramData|Public)') {
                $suspicious += "Value in user-writable path"
              }

              # Check 5: value path is an interpreter
              if ($value -match '(powershell|cmd.exe|wscript|cscript|mshta|rundll32|regsvr32)') {
                $suspicious += "Value invokes interpreter"
              }

              if ($suspicious.Count -gt 0) {
                [PSCustomObject]@{
                  Key   = $key
                  Name  = $name
                  Value = $value
                  Suspicious = ($suspicious -join '; ')
                }
              }
            }
          }
        }
    registry: |-
        Adversary value name patterns documented in
        real-world incidents:

        Spoofed Microsoft/system product names:
        - "Microsoft Defender", "WindowsDefender",
          "WinDefender" (often with path NOT in C:\Program
          Files\Windows Defender\)
        - "Windows Update", "WindowsUpdater", "WinUpdate"
        - "Microsoft Office Update Service"
        - "Microsoft Edge Update"
        - "OneDriveStandaloneUpdater"

        Spoofed third-party product names:
        - "Adobe Updater", "AdobeARMservice"
        - "Google Update", "GoogleUpdaterTaskMachine"
        - "iTunes Helper", "Apple Push"
        - "Skype Helper"
        - "Chrome Helper", "ChromeUpdater"

        Random / machine-generated names:
        - 8-16 character lowercase alphanumeric strings
        - GUID-format strings (sometimes)
        - Names matching the executable filename
          (poor opsec - the value name matching the EXE
          name is a tell)

        Blank/whitespace abuse:
        - Some malware uses empty string value names
        - Older Windows versions render these as blanks
          in regedit, making them visually 'hidden'
        - This is detectable; regedit lists them but
          with empty Name field

        Hunt with: combine the PowerShell scan with
        manual investigation of any flagged entries.
        Pattern-match alone has false positives -
        investigation confirms.
    ossdetect: |-
        Sigma:
        - registry_event_run_key_susp_value_name.yml
        - registry_event_run_key_susp_path.yml
        - registry_event_susp_run_key_spoofed_name.yml

        Atomic Red Team:
        - T1547.001 (tests cover varied value naming)

        Hayabusa:
        - RegistryRunKeyMod rules
        - Generic patterns; manual investigation needed
          for the spoofing-detection cases

        Velociraptor:
        - Windows.Sys.AutoRuns
        - Output includes Publisher field; cross-reference
          against ExpectedPublisher for each value name

        YARA / static analysis:
        - Many YARA rules exist for malware families
          that use specific value name patterns
        - Useful for tying a Run key entry back to a
          specific known family

        Reference: SANS DFIR
        'Windows Registry Forensics' poster
        - Comprehensive coverage of Run key abuse patterns
    tools: |-
        Adversary tradecraft using these patterns:

        Manual operators:
        - Choose value names to blend in with installed
          software fingerprint of the target
        - Adjust naming based on what runs on the host
          (recon first, persist with matching names)

        Commodity malware default names (don't always bother
        to spoof - just use whatever the loader picks):
        - Emotet variants: random or product-mimicking
        - IcedID: often impersonates Microsoft/Adobe
        - QakBot: random alphanumeric
        - Agent Tesla: varies by sample

        Red team tools:
        - SharPersist supports custom value names
        - Cobalt Strike aggressor scripts allow naming
          control for value entries
        - Empire persistence modules

        Legitimate software that this hunt may flag
        (tune accordingly):
        - Hardware vendor utilities with random IDs in
          value names (Dell, HP, Lenovo)
        - Some game launchers
        - IT management agents with non-standard install
          paths and generic-sounding names
    actors:
        - Lazarus
        - APT41
        - APT35
    analyst_notes: 'Value name analysis is a behavioral fingerprinting approach to autostart hunting - it asks ''does this entry look like something a legitimate vendor would create, or does it look like an adversary''s attempt to blend in?'' The detection is necessarily fuzzy because legitimate software does occasionally use random-looking or generic value names, and adversaries with good opsec choose names that match the actual installed software on the target host. The strongest signals combine multiple flags: a Microsoft-sounding value name pointing to a binary in AppData is high-confidence malicious; a Microsoft-sounding value name pointing to a binary in C:\Program Files\Microsoft is probably legitimate. The ''name spoof + path mismatch'' pattern is the highest-value single check because legitimate Microsoft software is consistently installed under Program Files paths, never in user-writable directories. Pair this hunt with the previous indicator''s general autostart sweep - this indicator adds the spoofing/naming dimension that pure path-based detection might miss. False positives are inherent here - run this as a triage prioritizer, not a final alert.'
