# sigma/rules/host/persistence/t1543_002_systemd_timer_persistence_timer_as_cron_alternat.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1543.002 - Systemd Timer Persistence (.timer as cron alternative) (Timer unit file creation)
id: 9d3a36ec-e5b5-53b2-8928-88b583792d4a
status: experimental
description: 'Creation of a .timer unit file as an alternative to cron-based persistence, or a .service unit with a corresponding .timer that triggers execution on a schedule, particularly when the timer unit or its target service ExecStart points to a binary in a writable or non-package path Query block: Timer unit file creation.'
references:
    - MITRE ATT&CK T1543.002
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1543.002
    - attack.g0096
    - attack.g0010
logsource:
    category: file_event
    product: windows
detection:
    selection:
        file.path:
            - '*/etc/systemd/system/*.timer'
            - '*/usr/lib/systemd/system/*.timer'
        event.action:
            - created
            - modified
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1543.002
    technique_name: 'Create or Modify System Process: Systemd Service'
    os: linux
    query_block: 1 of 2
    original_query: |-
        file.path: (*/etc/systemd/system/*.timer OR */usr/lib/systemd/system/*.timer)
        AND event.action: ("created" OR "modified")
    sysmon: |-
        // Auditd rules for timer unit creation
        -w /etc/systemd/system/ -p wa -k systemd_unit_write
        -w /usr/lib/systemd/system/ -p wa -k systemd_unit_write

        // Sysmon for Linux EID 11 (FileCreate)
        TargetFilename matches: *.timer
          AND path contains: /etc/systemd/ OR /usr/lib/systemd/

        // systemctl enable of a timer unit
        Image=*/systemctl AND CommandLine=*enable* AND CommandLine=*.timer*
    powershell: |-
        # Hunt for systemd timer persistence
        echo "[*] === All .timer units ==="
        systemctl list-timers --all --no-pager

        echo ""
        echo "[*] === Non-package .timer files ==="
        find /etc/systemd/system /usr/lib/systemd/system -name '*.timer' -type f 2>/dev/null | while read t; do
          pkg=$(dpkg -S "$t" 2>/dev/null || rpm -qf "$t" 2>/dev/null)
          if [ -z "$pkg" ] || echo "$pkg" | grep -q 'not found|not owned'; then
            echo "  UNPACKAGED: $t"
            grep -E 'OnCalendar|OnBoot|OnActiveSec' "$t" 2>/dev/null
            unit=$(grep 'Unit=' "$t" 2>/dev/null | cut -d= -f2)
            [ -n "$unit" ] && echo "  -> triggers: $unit"
          fi
        done
    registry: |-
        No registry artifact (Linux technique).

        Timer units use OnCalendar, OnBootSec, OnActiveSec directives
        for scheduling. Unlike cron, timers have built-in randomized
        delay (AccuracySec), persistent across reboots (Persistent=true),
        and can trigger any service unit.
    ossdetect: |-
        Sigma:
        - lnx_auditd_systemd_unit_file_creation.yml
        Wazuh:
        - FIM alert on /etc/systemd/system/*.timer
    tools: |-
        osquery:
          SELECT * FROM systemd_units WHERE id LIKE '%.timer';
        Wazuh FIM on /etc/systemd/system/
        Auditd rules for systemd directory writes
    actors:
        - APT41
        - Turla
    analyst_notes: 'Systemd timers are increasingly used by adversaries as a cron alternative because they offer more flexible scheduling, randomized jitter (AccuracySec), and are less frequently audited than crontab entries. The detection approach mirrors the .service detection: verify timer units against the package manager (dpkg -S / rpm -qf), flag unpackaged timers, and inspect what service they trigger.'
