# sigma/rules/host/persistence/t1546_003_suspicious_consumer_patterns_and_trigger_analysi.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: 'T1546.003 - Suspicious Consumer Patterns and Trigger Analysis (Real-time: dangerous consumer types)'
id: 5399000a-25e1-52ea-9bf6-fb85d86aa739
status: experimental
description: 'WMI consumer with execution-capable type pointing to interpreter, or filter with persistence-intent trigger (system boot, logon, perfdata polling) Query block: Real-time: dangerous consumer types.'
references:
    - MITRE ATT&CK T1546.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1546.003
    - attack.g0016
    - attack.g0096
    - attack.g0010
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 20
        winlog.event_data.Type:
            - CommandLineEventConsumer
            - ActiveScriptEventConsumer
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1546.003
    technique_name: 'Event Triggered Execution: WMI Event Subscription'
    os: win
    query_block: 1 of 2
    original_query: |-
        winlog.event_id: 20
        AND winlog.event_data.Type: ("CommandLineEventConsumer" OR "ActiveScriptEventConsumer")
    sysmon: |-
        // Real-time Sysmon detection of suspicious patterns
        // during subscription creation (EID 19/20/21):

        EventID=20 (Consumer Created)
        Type=CommandLineEventConsumer OR ActiveScriptEventConsumer
        // These two consumer types execute code -
        // other consumer types (NTEventLogEventConsumer,
        // SMTPEventConsumer, LogFileEventConsumer)
        // don't execute attacker code directly

        EventID=19 (Filter Created)
        // Suspicious filter queries indicating persistence:
        Query matches (any of):
          *__InstanceCreationEvent* (process / file create triggers)
          *__InstanceModificationEvent* (state change triggers)
          *Win32_LocalTime* (time-based triggers)
          *Win32_PerfFormattedData_PerfOS_System*
            (uptime-based triggers - fires after N seconds boot)
          *Win32_LogonSession* (logon triggers)
    powershell: |-
        # Hunt for suspicious consumer + filter pairings on live host

        # 1. Inventory all execution-capable consumers
        $execConsumers = Get-WMIObject -Namespace root\subscription -Class __EventConsumer |
          Where-Object {
            $_.__CLASS -eq 'CommandLineEventConsumer' -or
            $_.__CLASS -eq 'ActiveScriptEventConsumer'
          }

        # 2. Pull all filters
        $filters = Get-WMIObject -Namespace root\subscription -Class __EventFilter

        # 3. Pull all bindings (the link between filters and consumers)
        $bindings = Get-WMIObject -Namespace root\subscription -Class __FilterToConsumerBinding

        # 4. Reconstruct armed subscriptions (filter + consumer + binding triple)
        foreach ($binding in $bindings) {
          # Binding properties reference objects via __RELPATH
          $filterPath = $binding.Filter
          $consumerPath = $binding.Consumer

          $filterName = $filterPath -replace '.*Name="([^"]+)".*','$1'
          $consumerName = $consumerPath -replace '.*Name="([^"]+)".*','$1'

          $filter = $filters | Where-Object { $_.Name -eq $filterName }
          $consumer = $execConsumers | Where-Object { $_.Name -eq $consumerName }

          if ($consumer) {
            # This is an execution-capable armed subscription
            [PSCustomObject]@{
              FilterName = $filterName
              FilterQuery = $filter.Query
              ConsumerName = $consumerName
              ConsumerClass = $consumer.__CLASS
              Action = if ($consumer.CommandLineTemplate) {
                         $consumer.CommandLineTemplate
                       } else {
                         $consumer.ScriptText
                       }
              TriggerType = switch -Regex ($filter.Query) {
                '__InstanceCreationEvent' { 'Object created (process/file)' }
                '__InstanceModificationEvent' { 'State change' }
                'Win32_LocalTime' { 'Time-based' }
                'Win32_PerfFormattedData' { 'Uptime / perfdata polling' }
                'Win32_LogonSession' { 'Logon event' }
                default { 'Other' }
              }
            }
          }
        }
    registry: |-
        No registry artifacts - all data is in the
        WMI repository binary files.

        Consumer type analysis - what each does:

        CommandLineEventConsumer:
        - Runs an arbitrary command via the Win32_Process
          Create method when the bound filter fires
        - CommandLineTemplate property = the command
        - Highest-risk consumer type
        - Example: cmd /c C:\evil.exe

        ActiveScriptEventConsumer:
        - Runs inline VBScript or JScript when the bound
          filter fires
        - ScriptText property = the inline code
        - ScriptFileName property = optional file path
          (less common; most adversary use is inline)
        - High-risk; no file artifact when inline

        NTEventLogEventConsumer:
        - Writes to Windows Event Log
        - Not directly attacker-executable
        - Lower-risk classification

        SMTPEventConsumer:
        - Sends email
        - Not directly attacker-executable
        - Could be used for data exfiltration but rare

        LogFileEventConsumer:
        - Writes to a text file
        - Not directly attacker-executable

        Filter query patterns and their meaning:

        __InstanceCreationEvent + Win32_Process:
        - Fires when ANY process is created
        - Adversary use: 'every time notepad.exe runs,
          also run my code'

        __InstanceModificationEvent + Win32_PerfFormattedData_PerfOS_System:
        - Fires when system uptime crosses N seconds
        - Adversary use: 'run my code 5 minutes after
          every boot' (classic boot-persistence trigger)

        __InstanceCreationEvent + __TimerEvent / Win32_LocalTime:
        - Time-based triggers
        - Adversary use: scheduled-task-like behavior
          without a scheduled task artifact
    ossdetect: |-
        Sigma:
        - sysmon_wmi_susp_consumer_type.yml (EID 20)
        - sysmon_wmi_susp_filter_query.yml (EID 19)
        - sysmon_wmi_persistence_binding.yml (EID 21)

        Atomic Red Team:
        - T1546.003 (multiple subscription persistence tests)
        - T1546.003 Test #1 (CommandLineEventConsumer + filter)
        - T1546.003 Test #5 (ActiveScriptEventConsumer)

        Hayabusa:
        - WMIDangerousConsumerType rules
        - WMISuspFilterQuery detection

        Velociraptor:
        - Windows.Forensics.WMIPersistence
          (returns parsed subscriptions with consumer
          type and filter query exposed for analysis)

        Foundational research:
        - 'WMI Offense, Defense, and Forensics' -
          William Ballenthin, Matt Graeber, Claudiu
          Teodorescu (Mandiant, 2014)
        - The single most cited reference for WMI
          persistence detection
        - Still accurate today; the technique
          hasn't fundamentally changed

        DEF CON / Black Hat talks:
        - Various WMI persistence talks 2014-2018
        - All build on the Mandiant 2014 foundation
    tools: |-
        Subscription persistence frameworks:

        PowerSploit / Empire:
        - New-UserPersistenceOption -WMI
        - Creates CommandLineEventConsumer with logon trigger
        - Documented adversary tradecraft

        WMI-Persistence (Matt Graeber / Justin Warner):
        - PowerShell module for WMI persistence creation
        - Multiple trigger types supported
        - Reference for understanding the technique

        SharPersist:
        - /t wmi flag adds WMI subscription
        - Standard red team persistence option

        Cobalt Strike:
        - Aggressor scripts for WMI persistence
        - Multiple trigger and consumer combinations

        Manual operator pattern:
        - Create filter using stable trigger
          (perf data polling or logon event)
        - Create CommandLineEventConsumer pointing
          to staged payload
        - Bind them together
        - Subscription persists until explicitly removed
        - No registry entry, no scheduled task,
          no startup folder file - invisible to most
          standard persistence checks

        Defensive note - subscription removal:
        Remove-WmiObject -Namespace root\subscription
          -Class __FilterToConsumerBinding
          -Filter "Filter='...'"
        (removes the binding; filter/consumer remain
         unless removed separately)
    actors:
        - APT29
        - APT41
        - Turla
        - FIN6
    analyst_notes: 'This indicator narrows the WMI subscription hunt to the highest-risk subset: consumers with execution capability (CommandLineEventConsumer, ActiveScriptEventConsumer) paired with filters that have persistence intent (boot triggers, logon triggers, periodic polling). The ''trigger type'' classification in the PowerShell script is genuinely useful for triage: a boot-uptime trigger (Win32_PerfFormattedData_PerfOS_System) paired with a CommandLineEventConsumer running powershell is the classic adversary persistence pattern documented across many incidents. The reverse - benign trigger types like ''__InstanceModificationEvent + a tracked legitimate object'' paired with NTEventLogEventConsumer - is more often a legitimate monitoring or AV subscription. Subscription persistence is genuinely hard to detect without targeted hunting because Sysmon''s WMI monitoring (EID 19/20/21) requires explicit config-level opt-in and many Sysmon deployments don''t include it. Velociraptor''s Windows.Forensics.WMIPersistence artifact is the single most valuable tool for fleet-wide hunting of this technique. The Mandiant 2014 paper remains the foundational reference - if you want to dig deeper into the technique, it''s the canonical read.'
