# sigma/rules/host/persistence/t1547_009_lnk_file_in_startup_folder_lnk_file_creation_in.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1547.009 - LNK File in Startup Folder (LNK file creation in Startup folders)
id: 25115d36-3ba8-5d1e-bace-152884e547fa
status: experimental
description: 'LNK shortcut written to per-user or all-users Startup folder - file launched at every logon Query block: LNK file creation in Startup folders.'
references:
    - MITRE ATT&CK T1547.009
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1547.009
    - attack.g0007
    - attack.g0096
    - attack.g0032
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 11
        file.extension: lnk
        file.path:
            - '*\Start Menu\Programs\Startup\*'
            - '*\Start Menu\Programs\StartUp\*'
    condition: selection
falsepositives:
    - legitimate vendor software does create LNKs in the all-users Startup folder during install - signature checking on the LNK target binary helps filter these out
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1547.009
    technique_name: 'Boot or Logon Autostart Execution: Shortcut Modification'
    os: win
    query_block: 1 of 3
    original_query: |-
        winlog.event_id: 11
        AND file.extension: "lnk"
        AND file.path: (*\Start Menu\Programs\Startup\* OR *\Start Menu\Programs\StartUp\*)
    sysmon: |-
        // Real-time detection - LNK file create in Startup folder:
        EventID=11
        TargetFilename matches:
          *\Microsoft\Windows\Start Menu\Programs\Startup\*.lnk
          OR *\Microsoft\Windows\Start Menu\Programs\StartUp\*.lnk

        // Also catch non-LNK files in Startup folder
        // (some adversaries drop .exe/.bat/.vbs directly):
        EventID=11
        TargetFilename matches:
          *\Start Menu\Programs\Startup\*
          OR *\Start Menu\Programs\StartUp\*

        // LNK execution at logon - explorer parent + LNK target:
        EventID=1
        ParentImage=*\explorer.exe
        CommandLine matches: suspicious interpreter or path
        // (the LNK launches whatever its target says)
    powershell: |-
        # Hunt for LNK files in Startup folders + parse their targets

        $startupPaths = @(
          "$env:ProgramData\Microsoft\Windows\Start Menu\Programs\StartUp",
          "$env:APPDATA\Microsoft\Windows\Start Menu\Programs\Startup"
        )

        # Some systems have additional per-user paths - enumerate all users:
        $additionalPaths = Get-ChildItem 'C:\Users' -Directory -EA SilentlyContinue | ForEach-Object {
          Join-Path $_.FullName 'AppData\Roaming\Microsoft\Windows\Start Menu\Programs\Startup'
        } | Where-Object { Test-Path $_ }

        $allPaths = $startupPaths + $additionalPaths | Sort-Object -Unique

        $shell = New-Object -ComObject WScript.Shell

        foreach ($path in $allPaths) {
          if (-not (Test-Path $path)) { continue }

          Get-ChildItem $path -File | ForEach-Object {
            $suspicious = @()

            if ($_.Extension -eq '.lnk') {
              # Parse LNK target/args
              try {
                $lnk = $shell.CreateShortcut($_.FullName)
                $target = $lnk.TargetPath
                $args = $lnk.Arguments

                # Check 1: LNK target is an interpreter
                if ($target -match '(powershell|cmd\.exe|wscript|cscript|mshta|rundll32|regsvr32)') {
                  $suspicious += "Target invokes interpreter"
                }

                # Check 2: LNK target in user-writable path
                if ($target -match '(AppData|Temp|ProgramData|\\Users\\Public)') {
                  $suspicious += "Target in user-writable path"
                }

                # Check 3: Arguments contain risk keywords
                if ($args -match '(ExecutionPolicy.*Bypass|EncodedCommand|WindowStyle.*Hidden|vbscript:|javascript:|http://|https://)') {
                  $suspicious += "Risky arguments"
                }

                # Check 4: Target unsigned (if it exists)
                if ($target -and (Test-Path $target)) {
                  $sig = Get-AuthenticodeSignature $target -EA SilentlyContinue
                  if ($sig.Status -ne 'Valid') {
                    $suspicious += "Target unsigned ($($sig.Status))"
                  }
                }

                [PSCustomObject]@{
                  Path = $_.FullName
                  Type = 'LNK'
                  Target = $target
                  Arguments = $args
                  Created = $_.CreationTime
                  Modified = $_.LastWriteTime
                  Suspicious = if ($suspicious) { $suspicious -join '; ' } else { '' }
                }
              } catch {
                Write-Warning "Failed to parse $($_.FullName): $_"
              }
            } else {
              # Non-LNK file in Startup folder (always worth flagging)
              $suspicious += "Non-LNK file in Startup folder (.$($_.Extension))"

              [PSCustomObject]@{
                Path = $_.FullName
                Type = $_.Extension
                Target = $_.FullName
                Arguments = ''
                Created = $_.CreationTime
                Modified = $_.LastWriteTime
                Suspicious = $suspicious -join '; '
              }
            }
          }
        } | Sort-Object Suspicious -Descending | Format-Table -AutoSize
    registry: |-
        Startup folder paths:

        All-users (machine-wide, requires admin to write):
        C:\ProgramData\Microsoft\Windows\Start Menu\
          Programs\StartUp\

        Per-user (any user can write to their own):
        C:\Users\<user>\AppData\Roaming\Microsoft\
          Windows\Start Menu\Programs\Startup\

        Note the casing: 'StartUp' (machine) vs 'Startup'
        (per-user) - Microsoft's inconsistency. File system
        is case-insensitive so detection patterns work
        either way but it's worth knowing.

        Files in these folders launch at logon - any
        extension that Windows knows how to execute:
        - .lnk (most common - shortcut to actual binary)
        - .exe (direct binary)
        - .bat / .cmd (batch script)
        - .vbs / .vbe (VBScript)
        - .js / .jse (JScript)
        - .ps1 (PowerShell - only if execution policy allows)
        - .hta (HTA - launches via mshta.exe)

        LNK file structure (binary format):
        - Inspect with LECmd (Eric Zimmerman) or
          PowerShell WScript.Shell COM object (shown above)
        - Key fields:
          - TargetPath: the binary actually launched
          - Arguments: command-line arguments
          - WorkingDirectory: where target runs from
          - IconLocation: icon path (often spoofed -
            document icon on a cmd-launching LNK)
          - Description: tooltip (often blank for adversary
            LNKs)
          - Hotkey: keyboard shortcut (rare for malicious)

        Registry mirror of Startup folder location:
        HKCU\Software\Microsoft\Windows\CurrentVersion\
          Explorer\Shell Folders\Startup
        HKCU\Software\Microsoft\Windows\CurrentVersion\
          Explorer\User Shell Folders\Startup
        - Default values point to the per-user Startup path
        - Adversaries occasionally REDIRECT these to a
          different folder they control - rare but documented
        - Worth checking these registry values match the
          expected default path

        File creation timestamps:
        - LNK files have their own MACB timestamps
        - LNK files also EMBED target file timestamps
          inside the binary structure
        - Useful for forensic timeline reconstruction
    ossdetect: |-
        Sigma:
        - file_event_win_startup_folder_persistence.yml
        - file_event_win_lnk_creation_in_startup.yml
        - proc_creation_win_susp_explorer_child.yml
          (catches the execution side at logon)

        Atomic Red Team:
        - T1547.001 (includes Startup folder tests)
        - T1547.009 Test #1 (LNK file creation in Startup)
        - T1547.009 Test #2 (modified existing LNK target)

        Hayabusa:
        - StartupFolderFileWrite rules (EID 11)
        - LnkLaunchingInterpreter detection (EID 1)

        Velociraptor:
        - Windows.Sys.AutoRuns
          (enumerates Startup folder contents alongside
          Run keys and other autostart locations)
        - Windows.Forensics.Lnk
          (LNK file parser - extracts target, arguments,
          working dir, and embedded timestamps)

        Sysinternals autoruns.exe:
        - Logon tab includes Startup folder entries
        - Highlights non-Microsoft signed entries
        - Excellent for manual IR triage

        LECmd by Eric Zimmerman:
        - Standalone LNK file parser
        - Extracts every field including hidden ones
        - Free tool, essential for LNK forensic work
        - Useful when handling LNKs collected from
          a forensic image rather than live host
    tools: |-
        LNK creation methods:

        Built-in Windows:
        - Right-click > New > Shortcut (GUI)
        - PowerShell WScript.Shell COM object (most common
          for adversary use):
            $s = New-Object -ComObject WScript.Shell
            $lnk = $s.CreateShortcut("$startupPath\Update.lnk")
            $lnk.TargetPath = "C:\Windows\System32\powershell.exe"
            $lnk.Arguments = "-ep bypass -w hidden -file C:\path\payload.ps1"
            $lnk.Save()

        Adversary frameworks:
        - SharPersist (/t startupfolder)
        - Empire / Starkiller persistence modules
        - Cobalt Strike persistence aggressor scripts
        - Most commodity malware uses this technique:
          Emotet, IcedID, QakBot, AgentTesla, RemcosRAT,
          NjRAT, DarkComet, and many more

        Common adversary LNK patterns:
        - LNK named to mimic legitimate startup items:
          "Microsoft Office Quick Start.lnk"
          "OneDrive.lnk", "Adobe Updater.lnk"
        - TargetPath = powershell.exe or wscript.exe
        - Arguments contain encoded/obfuscated payload
        - IconLocation set to mimic the spoofed app icon
        - Working directory = user's home folder (looks normal)

        Detection note:
        - Per-user startup is more common because it requires
          no admin privileges - any compromised user can
          persist their own session
        - All-users startup requires admin and indicates
          more capable / privileged compromise

        Non-LNK files in Startup folder:
        - Adversaries occasionally drop .exe, .bat, .vbs
          directly without using LNK
        - Slightly noisier (less common in legitimate use)
        - Same detection - any file in Startup folder warrants
          inspection
    actors:
        - APT28
        - APT41
        - Lazarus
    analyst_notes: 'Startup folder persistence is the filesystem complement to Run key persistence - same effect (run at logon), different mechanism. The technique is universal: virtually every commodity malware family uses either Run keys, Startup folder LNKs, or both. The detection focus is twofold: real-time file creation in Startup paths (Sysmon EID 11) and hunting established LNKs with suspicious targets. The PowerShell LNK parser in the script is genuinely useful - it reads every LNK in every Startup path, parses the COM-accessible properties (TargetPath, Arguments, WorkingDirectory), and flags ones that point to interpreters or user-writable paths. Per-user Startup is more common than all-users because it requires no admin privileges - any compromised user account can persist itself. The non-LNK file check matters: while most adversary persistence uses LNK shortcuts, some malware drops .vbs or .bat files directly into the Startup folder. Any non-LNK file in a Startup folder is worth investigating because legitimate software almost always uses LNK shortcuts for startup entries. False positives: legitimate vendor software does create LNKs in the all-users Startup folder during install - signature checking on the LNK target binary helps filter these out.'
