# sigma/rules/host/persistence/t1547_001_hunt_established_autostart_persistence.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/persistence.js
title: T1547.001 - Hunt Established Autostart Persistence
id: e7ad1293-e5e7-5586-9464-67ecf09b5b68
status: experimental
description: Enumerate all known Run keys, RunOnce variants, Startup folders, Winlogon hooks, and IFEO debuggers - find what already persists
references:
    - MITRE ATT&CK T1547.001
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.persistence
    - attack.t1547.001
    - attack.g0016
    - attack.g0096
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        winlog.event_id:
            - 12
            - 13
        registry.path:
            - '*\Run\*'
            - '*\RunOnce\*'
            - '*\Userinit'
            - '*\Winlogon\Shell'
            - '*\Image File Execution Options\*'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/persistence.js
    domain: host
    technique: T1547.001
    technique_name: 'Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder'
    os: win
    original_query: |-
        winlog.event_id: (12 OR 13)
        AND registry.path: (*\Run\* OR *\RunOnce\* OR *\Userinit OR *\Winlogon\Shell OR *\Image File Execution Options\*)
    sysmon: |-
        // This indicator is hunting-oriented - sweeping
        // existing autostart locations rather than catching
        // new writes in real-time (previous indicator).
        //
        // Sysmon is not the primary tool here - PowerShell
        // enumeration and autoruns.exe are. However, EID 12
        // (registry key/value access) can supplement:

        EventID=12 (RegistryEvent Object create/delete)
        TargetObject=<any Run key path>
        // Useful for catching deletion events that hunting
        // might miss (adversary cleans up after staging)
    powershell: |-
        # Comprehensive autostart enumeration - live host
        # Covers Run keys, Startup folders, Winlogon, IFEO

        $results = @()

        # Run key variants (HKLM + HKCU + Wow6432Node)
        $runKeys = @(
          'HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\Run',
          'HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce',
          'HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\RunServices',
          'HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\RunServicesOnce',
          'HKLM:\SOFTWARE\Wow6432Node\Microsoft\Windows\CurrentVersion\Run',
          'HKLM:\SOFTWARE\Wow6432Node\Microsoft\Windows\CurrentVersion\RunOnce',
          'HKCU:\SOFTWARE\Microsoft\Windows\CurrentVersion\Run',
          'HKCU:\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce',
          'HKCU:\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Windows\Load',
          'HKCU:\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Windows\Run'
        )

        foreach ($key in $runKeys) {
          if (Test-Path $key) {
            $props = Get-ItemProperty $key
            $props.PSObject.Properties | Where-Object {
              $_.Name -notmatch '^PS' -and $_.Value
            } | ForEach-Object {
              $results += [PSCustomObject]@{
                Source = $key
                Name = $_.Name
                Value = $_.Value
              }
            }
          }
        }

        # Winlogon Userinit + Shell
        $winlogon = Get-ItemProperty 'HKLM:\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Winlogon'
        $results += [PSCustomObject]@{
          Source = 'Winlogon\Userinit'; Name = 'Userinit'; Value = $winlogon.Userinit
        }
        $results += [PSCustomObject]@{
          Source = 'Winlogon\Shell'; Name = 'Shell'; Value = $winlogon.Shell
        }

        # Image File Execution Options - any Debugger values
        Get-ChildItem 'HKLM:\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Image File Execution Options' |
          ForEach-Object {
            $dbg = (Get-ItemProperty $_.PSPath -Name Debugger -EA SilentlyContinue).Debugger
            if ($dbg) {
              $results += [PSCustomObject]@{
                Source = 'IFEO Debugger'
                Name = $_.PSChildName
                Value = $dbg
              }
            }
          }

        # Startup folder contents (both machine and user)
        $startupPaths = @(
          "$env:ProgramData\Microsoft\Windows\Start Menu\Programs\StartUp",
          "$env:APPDATA\Microsoft\Windows\Start Menu\Programs\Startup"
        )
        foreach ($path in $startupPaths) {
          if (Test-Path $path) {
            Get-ChildItem $path -File | ForEach-Object {
              $results += [PSCustomObject]@{
                Source = "StartupFolder: $path"
                Name = $_.Name
                Value = $_.FullName
              }
            }
          }
        }

        # Output and flag suspicious
        $results | ForEach-Object {
          $_ | Add-Member -NotePropertyName Suspicious -NotePropertyValue (
            $_.Value -match '(AppData|Temp|ProgramData|\\Users\\Public|powershell|cmd.exe|wscript|cscript|mshta|rundll32|regsvr32)'
          ) -PassThru
        } | Sort-Object Suspicious -Descending | Format-Table -AutoSize
    registry: |-
        Same registry paths as previous indicator -
        the difference is approach: real-time vs hunting.

        Hunting workflow:
        1. Run the PowerShell sweep above on the target host
        2. Diff against a known-good baseline if available
        3. For each result, check:
           - Is the value path signed? (Authenticode check)
           - Is the value path in a standard install location?
           - Does the value name match known legitimate
             software, or is it spoofed/random?
           - When was the registry key last modified?
             (Get-Item <key> | Select LastWriteTime)
        4. Suspicious entries: pivot to file analysis
           of the target binary

        Baseline approach:
        - Export current autostart inventory from a clean
          reference machine: autoruns.exe -a * -ct
        - Save as baseline CSV
        - Compare against suspect systems with:
          Compare-Object (Import-CSV baseline.csv)
            (Import-CSV suspect.csv) -Property Image
        - Diff highlights additions and deletions

        Last-write timestamps:
        - Each registry KEY has a LastWriteTime
        - Useful for narrowing 'when was this persistence
          installed' during IR triage
        - Available via:
          (Get-Item <path>).LastWriteTime
          or RegRipper / RECmd for offline registry analysis
    ossdetect: |-
        Sigma:
        - Detection rules in Sigma cover real-time writes
          (previous indicator) - hunting is a manual sweep
          process, not a real-time rule

        Atomic Red Team:
        - T1547.001 (multiple tests reproduce persistence)
        - Use to validate the hunt: install via Atomic,
          then run the PowerShell sweep to confirm it surfaces

        Hayabusa:
        - Hayabusa rules cover the real-time side
        - For hunting use, run autoruns/Velociraptor

        Velociraptor:
        - Windows.Sys.AutoRuns (dedicated artifact)
        - The gold standard for this type of sweep
        - Returns parsed autostart entries with
          signature and reputation context

        Open-source autoruns hunters:
        - Various community PowerShell scripts on GitHub
        - ThreatHunting repos by F-Secure, MITRE, Splunk

        Reference: SANS poster
        'Hunt Evil: Persistence Mechanisms'
        - Updated annually with new persistence locations
        - Good single-page reference for the technique surface
    tools: |-
        Sysinternals autoruns.exe:
        - The definitive tool for this kind of sweep
        - Run with: autoruns.exe -a * -accepteula -nobanner -ct > out.csv
        - Output columns: Entry Location, Entry, Enabled,
          Description, Publisher, Image Path, Version,
          Time, Hash (MD5/SHA1/SHA256), Virus Total
        - Flags unsigned entries in red
        - Can submit hashes directly to VirusTotal
        - For live host triage - close to mandatory

        Velociraptor (fleet-wide):
        - Windows.Sys.AutoRuns
        - Runs autoruns equivalent across many hosts
        - Best for environment-wide persistence sweeps
        - Output normalized for cross-host comparison

        KAPE (Eric Zimmerman):
        - Targets/Compound/AutomaticDestinations.tkape
        - AutomaticDestinations target captures Recent
          files / jump lists which can reveal persistence
        - Useful for IR forensic timeline building

        RegRipper (offline registry analysis):
        - Parses registry hives extracted from images
        - Plugins for Run keys, Winlogon, services, etc.
        - Useful when working from forensic image vs live host

        PowerShell native (the script above):
        - No tool deploy required
        - Works on any modern Windows
        - Good for ad-hoc remote IR via WinRM
    actors:
        - APT29
        - APT41
    analyst_notes: 'Hunting established autostart persistence is one of the most fundamental Windows IR skills - if a host is suspected of being compromised, this sweep is almost always the first or second step. Sysinternals autoruns.exe is the canonical tool because Microsoft maintains its list of autostart locations and it covers 200+ extension points that PowerShell scripts miss. The PowerShell script above covers the high-value subset and is useful when you can''t or shouldn''t deploy autoruns.exe to a host (forensic preservation, locked-down environments, remote IR via WinRM where pushing a binary is awkward). The baseline-and-diff approach is the highest-fidelity hunting method: capture autoruns output on a known-good reference image, then diff against suspect hosts to surface only the entries that don''t appear in the baseline. This collapses a 200-entry autoruns output into 5-10 entries worth investigating. The LastWriteTime trick for registry keys is genuinely useful in incident timelining - it tells you when the persistence was established, which often correlates with initial compromise or lateral movement events visible in other logs.'
