# sigma/rules/host/execution/t1047_wmi_event_subscription_persistence_alert_on_any.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/execution.js
title: T1047 - WMI Event Subscription Persistence (Alert on any WMI subscription activity (EID 19/20/21))
id: 25a1575a-3a68-5623-a82d-5c57dd82a2f8
status: experimental
description: 'WMI EventFilter, EventConsumer, or FilterToConsumerBinding creation - fileless persistence via WMI repository Query block: Alert on any WMI subscription activity (EID 19/20/21).'
references:
    - MITRE ATT&CK T1047
    - T1546.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.execution
    - attack.t1047
    - attack.g0016
    - attack.g0096
    - attack.g0010
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id:
            - 19
            - 20
            - 21
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/execution.js
    domain: host
    technique: T1047
    technique_name: Windows Management Instrumentation
    os: win
    query_block: 1 of 3
    original_query: 'winlog.event_id: (19 OR 20 OR 21)'
    sysmon: |-
        // Sysmon has three dedicated WMI persistence event IDs:

        EventID=19 (WmiEventFilter activity detected)
        // Fires when a WMI event filter is registered
        // Filter = the trigger condition (e.g., system boot,
        // process creation, time interval)
        // Key fields: Name, Query, QueryLanguage

        EventID=20 (WmiEventConsumer activity detected)
        // Fires when a WMI event consumer is registered
        // Consumer = the action to take when filter fires
        // Two dangerous consumer types:
        //   CommandLineEventConsumer (runs a command)
        //   ActiveScriptEventConsumer (runs VBScript/JScript)
        // Key fields: Name, Type, Destination/ScriptText

        EventID=21 (WmiEventConsumerToFilter activity)
        // Fires when filter and consumer are bound together
        // This is the final step that arms the subscription
        // Key fields: Consumer, Filter
    powershell: |-
        # Hunt for WMI event subscriptions via Sysmon EID 19/20/21
        19, 20, 21 | ForEach-Object {
          $id = $_
          Get-WinEvent -FilterHashtable @{
            LogName='Microsoft-Windows-Sysmon/Operational';
            ID=$id
          } -ErrorAction SilentlyContinue
        } | Select TimeCreated, Id,
          @{n='EventType';e={
            switch ($_.Id) {
              19 { 'WMI Filter Registered' }
              20 { 'WMI Consumer Registered' }
              21 { 'Filter-Consumer Binding' }
            }
          }},
          @{n='Details';e={$_.Message}} |
          Sort-Object TimeCreated -Descending

        # Direct WMI query for active subscriptions
        # (run on live host or via Velociraptor/remote PS)
        Get-WMIObject -Namespace rootsubscription -Class __EventFilter |
          Select Name, Query, QueryLanguage

        Get-WMIObject -Namespace rootsubscription -Class __EventConsumer |
          Select Name, __CLASS, CommandLineTemplate, ScriptText

        Get-WMIObject -Namespace rootsubscription -Class __FilterToConsumerBinding |
          Select Filter, Consumer
    registry: |-
        WMI persistence does NOT use registry run keys or
        scheduled task XML - it lives in the WMI repository:

        C:\Windows\System32\wbem\Repository\
          OBJECTS.DATA  - main object store
          INDEX.BTR     - index file
          INDEX.MAP     - index map

        Not directly human-readable - use these tools:
        - PyWMIPersistenceFinder (David Reaves / mandiant)
          python PyWMIPersistenceFinder.py OBJECTS.DATA
        - WMIForensics (python, similar capability)
        - Velociraptor artifact Windows.Forensics.WMIPersistence
        - autoruns.exe (Sysinternals) - lists WMI subscriptions
          under the WMI tab

        Legitimate WMI subscriptions exist in healthy
        environments (SCCM, monitoring agents, AV products)
        - baseline before alerting on all subscriptions.
        Known-good namespaces: rootccm, rootcimv2sms
        Adversary subscriptions typically use: rootsubscription
    ossdetect: |-
        Sigma:
        - sysmon_wmi_event_subscription.yml
        - sysmon_wmi_persistence_script_event_consumer.yml
        - proc_creation_win_wmic_eventconsumer_creation.yml

        Atomic Red Team:
        - T1546.003 (WMI event subscription tests)
        - T1047 (WMI execution tests)

        Hayabusa:
        - WMIEventSubscription (Sysmon EID 19/20/21) rules

        Velociraptor:
        - Windows.Persistence.PermanentWMIEvents
        - Windows.Detection.WMIProcessCreation

        Sysinternals autoruns.exe:
        - WMI tab enumerates __EventFilter / __EventConsumer /
          __FilterToConsumerBinding entries
    tools: |-
        PowerSploit / PowerShell Empire:
        - New-UserPersistenceOption -WMI
        - Installs CommandLineEventConsumer triggering on logon

        Metasploit:
        - exploit/windows/local/wmi_persistence

        Sharp-WMI (C# WMI execution tool)
        WMImplant (PowerShell WMI C2 framework)
        Cobalt Strike (WMI persistence module)

        Manual WMI subscription (PowerShell):
        $filter = Set-WmiInstance -Namespace rootsubscription
          -Class __EventFilter -Arguments @{
            Name='Updater';
            QueryLanguage='WQL';
            Query='SELECT * FROM __InstanceModificationEvent
              WITHIN 60 WHERE TargetInstance ISA "Win32_PerfFormattedData_PerfOS_System"'
          }
        $consumer = Set-WmiInstance -Namespace rootsubscription
          -Class CommandLineEventConsumer -Arguments @{
            Name='Updater';
            CommandLineTemplate='cmd /c payload.exe'
          }
        Set-WmiInstance -Namespace rootsubscription
          -Class __FilterToConsumerBinding -Arguments @{
            Filter=$filter; Consumer=$consumer
          }
    actors:
        - APT29
        - APT41
        - FIN6
        - Turla
    analyst_notes: 'WMI event subscription persistence is the most forensically evasive standard persistence technique on Windows. There is no registry run key, no scheduled task XML, no startup folder file - the persistence object lives entirely inside the WMI repository binary (OBJECTS.DATA), which most IR tools and AV products do not inspect by default. Sysmon EID 19/20/21 are the primary detection mechanism and are only logged if Sysmon is deployed with WMI monitoring enabled in its config - verify your Sysmon config actually captures these before assuming you have coverage. The three-event chain (19 = filter registered, 20 = consumer registered, 21 = binding created) tells the complete story: what triggers it, what it does, and that it is now armed. In practice, EID 21 (binding) is the highest-value single alert because a filter and consumer registered independently are less meaningful than when they are bound together. Two consumer types are dangerous: CommandLineEventConsumer (runs a shell command) and ActiveScriptEventConsumer (runs VBScript or JScript inline - no file on disk at all). Legitimate WMI subscriptions exist (SCCM, some AV products) - baseline your environment before alerting on all subscription activity. The root\subscription namespace is the canonical adversary location; legitimate subscriptions more often live under root\cimv2 or vendor-specific namespaces.'
