# sigma/rules/host/execution/t1053_003_cron_job_crontab_abuse_cron_file_modification.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/execution.js
title: T1053.003 - Cron Job / Crontab Abuse (Cron file modification)
id: b720931b-e1d3-56b1-a150-29ba61792c5a
status: experimental
description: 'Creation or modification of a crontab or /etc/cron.* entry invoking a shell/interpreter or a payload in a writable path - or cron spawning a suspicious child Query block: Cron file modification.'
references:
    - MITRE ATT&CK T1053.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.execution
    - attack.t1053.003
    - attack.g0096
    - attack.g0007
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 11
        file.path:
            - /etc/crontab
            - /etc/cron.d/*
            - /etc/cron.hourly/*
            - /etc/cron.daily/*
            - /etc/cron.weekly/*
            - /var/spool/cron/*
            - /etc/systemd/system/*.timer
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/execution.js
    domain: host
    technique: T1053.003
    technique_name: 'Scheduled Task/Job: Cron'
    os: linux
    query_block: 1 of 3
    original_query: |-
        winlog.event_id: 11
        AND file.path: ("/etc/crontab" OR "/etc/cron.d/*" OR "/etc/cron.hourly/*" OR "/etc/cron.daily/*" OR "/etc/cron.weekly/*" OR "/var/spool/cron/*" OR "/etc/systemd/system/*.timer")
    sysmon: |-
        // Sysmon for Linux: file modify (EID 11) of cron locations,
        // and ProcessCreate (EID 1) of cron-spawned children.
        EventID=11 (FileCreate / modify)
        TargetFilename matches:
          /etc/crontab
          /etc/cron.d/*  /etc/cron.hourly/*  /etc/cron.daily/*
          /var/spool/cron/*  /var/spool/cron/crontabs/*
          /etc/systemd/system/*.timer

        EventID=1 (cron-spawned execution)
        ParentImage=*/cron OR */crond OR */CRON
        Image=*/bash OR */sh OR */python* OR */curl OR */wget
        CommandLine = a payload / download / reverse shell
    powershell: |-
        # (Auditd / Shell hunt - Linux row)

        # Enumerate ALL cron entries across the system (review each)
        for u in $(cut -f1 -d: /etc/passwd); do
          crontab -l -u "$u" 2>/dev/null | grep -v '^#' | sed "s/^/[$u] /"
        done
        cat /etc/crontab 2>/dev/null
        ls -la /etc/cron.d/ /etc/cron.hourly/ /etc/cron.daily/ \
          /etc/cron.weekly/ /etc/cron.monthly/ 2>/dev/null
        cat /etc/cron.d/* 2>/dev/null
        ls -la /var/spool/cron/ /var/spool/cron/crontabs/ 2>/dev/null

        # systemd timers (the modern cron)
        systemctl list-timers --all
        for t in /etc/systemd/system/*.timer; do echo "== $t =="; cat "$t"; done 2>/dev/null

        # auditd watch rules for cron tampering:
        #   -w /etc/crontab -p wa -k cron_persist
        #   -w /etc/cron.d/ -p wa -k cron_persist
        #   -w /var/spool/cron/ -p wa -k cron_persist
        ausearch -k cron_persist -i 2>/dev/null

        # Flag cron entries invoking shells/downloads/temp paths
        grep -RnE 'curl|wget|base64|/dev/tcp/|/tmp/|/dev/shm/|bash -i' \
          /etc/crontab /etc/cron.* /var/spool/cron 2>/dev/null
    registry: |-
        (File Artifacts - Linux row)

        No registry. Cron artifact locations (all worth a sweep):

        System-wide:
          /etc/crontab               (master system crontab)
          /etc/cron.d/*              (drop-in system cron files)
          /etc/cron.hourly|daily|weekly|monthly/*  (script dirs)

        Per-user:
          /var/spool/cron/<user>            (RHEL/CentOS)
          /var/spool/cron/crontabs/<user>   (Debian/Ubuntu)

        systemd timers (modern equivalent):
          /etc/systemd/system/*.timer + matching *.service
          /usr/lib/systemd/system/*.timer
          ~/.config/systemd/user/*.timer  (user timers)

        Logs:
          /var/log/cron , /var/log/syslog (cron execution lines)
          journalctl -u cron / journalctl -u crond

        Investigation pivots:
        - Any cron entry invoking curl/wget/base64, a reverse
          shell, or a payload in /tmp,/dev/shm is the signature
        - @reboot entries (run at boot) are favored for persistence
        - Dot-prefixed or odd-named files in /etc/cron.d
        - A user crontab for an account that shouldn't have one
          (service accounts, nologin users)
        - systemd .timer pointing at a .service that runs a
          writable-path binary
    ossdetect: |-
        Sigma (Linux rules):
        - lnx_auditd_cron_file_modification.yml
        - lnx_cron_susp_command.yml
        - lnx_systemd_timer_creation.yml

        Atomic Red Team:
        - T1053.003 (cron job tests)

        Auditd:
        - -w watches on /etc/crontab, /etc/cron.d/, /var/spool/cron

        Sysmon for Linux:
        - FileCreate (EID 11) on cron paths; ProcessCreate (EID 1)
          for cron-spawned children

        Falco:
        - "Schedule Cron Jobs" / write-below-etc rules

        Velociraptor:
        - Linux.Sys.Crontab artifact (fleet-wide cron enumeration)
    tools: |-
        Built-in: crontab -e / direct file writes (no tooling
          needed - cron abuse is a few lines)
        Metasploit - cron persistence modules
        Empire / Sliver Linux persistence modules
        LinPEAS / linux-smart-enumeration - flag writable cron
          files and misconfigured cron (recon)
        GTFOBins - entries abusing cron-invoked SUID helpers
        Manual operators - @reboot cron and /etc/cron.d drops are
          among the most common Linux persistence/exec methods
    actors:
        - APT41
        - APT28
        - TeamTNT
        - Kinsing
        - Rocke
    analyst_notes: 'Cron is both an execution and a persistence technique; this row focuses on detecting the abuse, and it pairs with the eventual Linux persistence coverage. The detection surface is broad because cron has many homes - the master /etc/crontab, drop-in files in /etc/cron.d, the script directories (/etc/cron.hourly|daily|weekly|monthly), per-user spools (/var/spool/cron on RHEL, /var/spool/cron/crontabs on Debian), and increasingly systemd .timer units which are the modern replacement and easy to overlook. Two detection angles: file modification on any of those locations (auditd watch rules keyed for easy searching, or Sysmon for Linux EID 11), and cron-spawned children that are suspicious (cron/crond as parent of a shell, python, curl, or a /tmp payload). The content tells are the same as elsewhere - curl/wget/base64, reverse shells, payloads in writable paths - plus the cron-specific @reboot directive favored for boot persistence. The proactive play is a full cron enumeration sweep across all users and all locations including systemd timers, flagging any entry that invokes a downloader, a shell, or a writable-path binary, and any crontab belonging to a service/nologin account that shouldn''t have one. This needs auditd watch rules on the cron paths to catch reliably in real time - note that as a prerequisite. Baseline legitimate scheduled maintenance (logrotate, package-update timers, monitoring agents) before alerting.'
