# sigma/rules/host/execution/t1053_003_system_cron_d_injection_and_spoofed_name_cron_en_1.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/execution.js
title: T1053.003 - System cron.d Injection and Spoofed-Name Cron Entries (crontab command execution)
id: dbf6872d-28a4-5227-b197-54af46d58389
status: experimental
description: 'A new file in /etc/cron.d/ or modification of /etc/crontab with an entry using a spoofed system-sounding name (apache, sshd, sysstat) or an entry whose command includes curl/wget piped to shell, base64 decoding, or targets a binary in /tmp or /dev/shm Query block: crontab command execution.'
references:
    - MITRE ATT&CK T1053.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.execution
    - attack.t1053.003
    - attack.g0096
    - attack.g1017
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        process.name: crontab
        process.args:
            - -l
            - -e
            - -r
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/execution.js
    domain: host
    technique: T1053.003
    technique_name: 'Scheduled Task/Job: Cron'
    os: linux
    query_block: 2 of 2
    original_query: |-
        process.name: "crontab"
        AND process.args: ("-l" OR "-e" OR "-r")
    sysmon: |-
        // Auditd rules for system cron directories
        -w /etc/cron.d/ -p wa -k cron_persist
        -w /etc/crontab -p wa -k cron_persist
        -w /var/spool/cron/ -p wa -k cron_persist

        // Sysmon for Linux EID 11 - file creation in cron directories
        TargetFilename matches: /etc/cron.d/*
          OR /etc/cron.daily/* OR /etc/cron.hourly/*
    powershell: |-
        # Hunt for suspicious system cron entries
        echo "[*] === /etc/cron.d contents ==="
        ls -la /etc/cron.d/

        echo ""
        echo "[*] === Non-package cron.d files ==="
        for f in /etc/cron.d/*; do
          [ -f "$f" ] || continue
          pkg=$(dpkg -S "$f" 2>/dev/null || rpm -qf "$f" 2>/dev/null)
          if [ -z "$pkg" ] || echo "$pkg" | grep -q 'not found\|not owned'; then
            echo "  UNPACKAGED: $f"
            cat "$f" | grep -v '^#\|^$'
            echo "  ---"
          fi
        done

        echo ""
        echo "[*] === Suspicious cron commands ==="
        grep -rh 'curl\|wget\|base64\|/dev/shm\|/tmp/' /etc/cron* /var/spool/cron/ 2>/dev/null | grep -v '^#'
    registry: |-
        No registry artifact (Linux technique).

        System cron hierarchy:
        - /etc/crontab: master system crontab
        - /etc/cron.d/: drop-in cron files (package-managed and custom)
        - /etc/cron.{hourly,daily,weekly,monthly}/: run-parts dirs
        - /var/spool/cron/crontabs/: per-user crontabs
    ossdetect: |-
        Sigma:
        - lnx_auditd_cron_job_creation.yml
        osquery:
        - crontab table sweep
        Wazuh:
        - FIM on /etc/cron.d/
    tools: |-
        osquery:
          SELECT * FROM crontab;
        Wazuh FIM on /etc/cron.d/
        Auditd rules for cron directory writes
    actors:
        - APT41
        - Volt Typhoon
    analyst_notes: 'System cron.d injection is preferred over user crontab by sophisticated operators because (a) /etc/cron.d files can specify the run-as user (root), (b) they survive user account changes, and (c) they are less visible than crontab -l for the affected user. Spoofed names (apache, sshd, sysstat) are a deliberate masquerade technique. The package-verification approach works: enumerate /etc/cron.d files, verify against the package manager, flag unpackaged entries.'
