# sigma/rules/host/execution/t1053_005_remote_scheduled_task_creation_source_schtasks_t.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/execution.js
title: 'T1053.005 - Remote Scheduled Task Creation (Source: schtasks targeting remote host)'
id: a55de688-676b-5270-9230-5c43f3871295
status: experimental
description: 'schtasks.exe /create with /s flag targeting remote host - lateral movement or remote persistence via Task Scheduler Query block: Source: schtasks targeting remote host.'
references:
    - MITRE ATT&CK T1053.005
    - T1021
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.execution
    - attack.t1053.005
    - attack.g0016
    - attack.g0096
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 1
        process.name: schtasks.exe
        process.command_line|all:
            - '*/create*'
            - '*/s *'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: high
tonk:
    fidelity: full
    source_file: host/js/data/execution.js
    domain: host
    technique: T1053.005
    technique_name: 'Scheduled Task/Job: Scheduled Task'
    os: win
    query_block: 1 of 3
    original_query: |-
        winlog.event_id: 1
        AND process.name: "schtasks.exe"
        AND process.command_line: (*/create* AND */s *)
    sysmon: |-
        // On SOURCE host:
        EventID=1
        Image=*\schtasks.exe
        CommandLine=*/create*
        AND CommandLine=*/s *
        // /s <hostname/IP> = remote target

        // Network connection from source:
        EventID=3
        Image=*\schtasks.exe
        DestinationPort=445
        // (schtasks remote uses SMB - TCP 445)

        // On DESTINATION host - Security Event:
        EventID=4698 (task created)
        // SubjectUserName will be the remote authenticating user
        // SubjectLogonId correlates to a Type 3 network logon
        // (Security EID 4624 LogonType=3 around same time)
    powershell: |-
        # Hunt for remote schtasks /create (source host)
        Get-WinEvent -FilterHashtable @{
          LogName='Microsoft-Windows-Sysmon/Operational';
          ID=1
        } | Where-Object {
          $_.Properties[4].Value -like '*\schtasks.exe' -and
          $_.Properties[10].Value -match '/create' -and
          $_.Properties[10].Value -match '/ss+S+'
        } | Select TimeCreated,
          @{n='CmdLine';e={$_.Properties[10].Value}},
          @{n='User';e={$_.Properties[12].Value}},
          @{n='Parent';e={($_.Properties[20].Value -split '\\')[-1]}} |
          Sort-Object TimeCreated -Descending

        # SMB connections from schtasks (EID 3)
        Get-WinEvent -FilterHashtable @{
          LogName='Microsoft-Windows-Sysmon/Operational';
          ID=3
        } | Where-Object {
          $_.Properties[4].Value -like '*\schtasks.exe' -and
          $_.Properties[16].Value -eq '445'
        } | Select TimeCreated,
          @{n='DestIP';e={$_.Properties[14].Value}},
          @{n='User';e={$_.Properties[12].Value}}
    registry: |-
        On DESTINATION host - task XML artifact:
        - C:\Windows\System32\Tasks\<TaskName>
          Created by the remote scheduling operation
          SubjectUserName in EID 4698 is the remote user -
          correlate with Security EID 4624 (LogonType=3)
          to confirm which account was used and from where

        Lateral movement context:
        - Remote scheduled task creation requires:
          1. Valid credentials on the target (or pass-the-hash)
          2. SMB access to target (TCP 445)
          3. Task Scheduler service running on target
        - This is the same credential/access requirement
          as PsExec, but uses a different protocol path
          and leaves a different artifact set

        Network pivot:
        - SMB connection (TCP 445) from source to target
          is visible in Zeek conn.log and Suricata alerts
        - Cross-reference with /net/lateral.html (Lateral Movement)
          reference for the network-side complement
    ossdetect: |-
        Sigma:
        - proc_creation_win_schtasks_remote.yml
        - win_security_scheduled_task_creation.yml
          (EID 4698 with network logon correlation)
        - network_connection_win_schtasks_smb.yml

        Atomic Red Team:
        - T1053.005 Test #4 (remote scheduled task)

        Hayabusa:
        - RemoteScheduledTask rules
        - SMB lateral movement correlation rules

        Velociraptor:
        - Windows.EventLogs.Sysmon (both hosts)
        - Windows.System.ScheduledTasks (destination)
        - Windows.Forensics.Autoruns (destination)
    tools: |-
        schtasks.exe /create /s <target> (built-in)
        Impacket atexec.py
          - Implements remote task scheduling over SMB
          - Does not use schtasks.exe on source host
          - Creates task, runs it, deletes it (less persistent)
        CrackMapExec --exec-method atexec
        PowerShell Register-ScheduledTask with
          -CimSession (remote CIM session)
        Cobalt Strike (remote task scheduling module)

        Key distinction: Impacket atexec and CrackMapExec
        atexec do NOT generate schtasks.exe EID 1 on source -
        detection shifts entirely to destination EID 4698
        and the SMB/authentication artifacts.

        atexec pattern: creates task, runs it immediately,
        deletes it - watch for short-lived tasks (EID 4698
        task created followed quickly by EID 4699 task deleted)
    actors:
        - APT29
        - APT41
    analyst_notes: 'Remote scheduled task creation is a lateral movement technique that predates most modern C2 frameworks - the Windows ''at'' command (now deprecated) and its successor schtasks /s have been used for decades. The detection split between source and destination is the same pattern as remote WMI: the source generates a process creation event (schtasks /s) and a network connection (SMB to 445), while the destination generates a task creation event (Security EID 4698) with a network logon (Security EID 4624 LogonType=3). The EID 4699 (task deleted) paired shortly after EID 4698 is worth a separate alert - it''s the signature of execution-only remote task use (Impacket atexec pattern) where the adversary creates a task, runs it immediately, and deletes it to minimize forensic evidence. A task that exists for less than 60 seconds and was created by a remote user is high-confidence malicious activity.'
