# sigma/rules/host/execution/t1053_005_scheduled_task_action_pointing_to_suspicious_bin.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/execution.js
title: 'T1053.005 - Scheduled Task Action Pointing to Suspicious Binary (Task execution: svchost (scheduler) spawning suspicious'
id: e64a9d75-f821-5a99-a4f1-8496e4a5d714
status: experimental
description: 'Existing or newly created scheduled task with action path in user-writable directory or running a scripting interpreter Query block: Task execution: svchost (scheduler) spawning suspicious child.'
references:
    - MITRE ATT&CK T1053.005
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.execution
    - attack.t1053.005
    - attack.g0016
    - attack.g0096
    - attack.g0032
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 1
        process.parent.name: svchost.exe
        process.parent.command_line: '*Schedule*'
        process.command_line:
            - '*\AppData\*'
            - '*\Temp\*'
            - '*\ProgramData\*'
            - '*powershell*'
            - '*wscript*'
            - '*mshta*'
            - '*rundll32*'
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/execution.js
    domain: host
    technique: T1053.005
    technique_name: 'Scheduled Task/Job: Scheduled Task'
    os: win
    query_block: 1 of 2
    original_query: |-
        winlog.event_id: 1
        AND process.parent.name: "svchost.exe"
        AND process.parent.command_line: *Schedule*
        AND process.command_line: (*\AppData\* OR *\Temp\* OR *\ProgramData\* OR *powershell* OR *wscript* OR *mshta* OR *rundll32*)
    sysmon: |-
        // Sysmon does not directly inspect task XML content -
        // use Windows Security Event 4698 for this.
        // Sysmon contribution: catch the EXECUTION of the task:

        EventID=1
        // Task runs as SYSTEM or as a specific user account
        // Parent will be taskeng.exe (legacy) or svchost.exe
        // (modern Task Scheduler host):
        ParentImage=*\svchost.exe
        ParentCommandLine=*Schedule*
        // AND child Image is suspicious:
        Image matches:
          *\AppData\* OR *\Temp\* OR *\ProgramData\*
          OR *\powershell.exe OR *\wscript.exe
          OR *\mshta.exe OR *\rundll32.exe
    powershell: |-
        # Enumerate all scheduled tasks and flag suspicious actions
        Get-ScheduledTask | ForEach-Object {
          $task = $_
          $actions = $task.Actions
          foreach ($action in $actions) {
            if ($action.Execute -match
              '(AppData|Temp|ProgramData|Public|powershell|wscript|cscript|mshta|rundll32|regsvr32|certutil)') {
              [PSCustomObject]@{
                TaskName    = $task.TaskName
                TaskPath    = $task.TaskPath
                Execute     = $action.Execute
                Arguments   = $action.Arguments
                State       = $task.State
                Author      = $task.Author
                RunAs       = $task.Principal.UserId
                LastRun     = $task.LastRunTime
                NextRun     = $task.NextRunTime
              }
            }
          }
        } | Sort-Object TaskPath

        # Also check raw XML for tasks that Get-ScheduledTask
        # may not fully surface (e.g., broken/hidden tasks):
        Get-ChildItem C:\Windows\System32\Tasks -Recurse -File |
          ForEach-Object {
            $content = Get-Content $_.FullName -Raw
            if ($content -match
              '(AppData|Temp|ProgramData|powershell|wscript|mshta|rundll32)') {
              [PSCustomObject]@{
                File    = $_.FullName
                Snippet = ($content -split '
        ' |
                  Select-String 'AppData|Temp|ProgramData|powershell|wscript|mshta|rundll32' |
                  Select-Object -First 3) -join '; '
              }
            }
          }
    registry: |-
        Task XML files - primary forensic artifact:
        - C:\Windows\System32\Tasks\<TaskName>
          Read with Get-Content or any text editor
          Key XML elements to inspect:
          <Exec><Command> - the binary being run
          <Exec><Arguments> - command line arguments
          <Principal><UserId> - run-as account
          <Triggers> - what causes the task to fire
          <RegistrationInfo><Author> - who created it

        Registry task cache (secondary):
        - HKLM\SOFTWARE\Microsoft\Windows NT\
          CurrentVersion\Schedule\TaskCache\Tasks\{GUID}\
          Actions value contains the task action in binary
          format - use specialized tools to parse

        Deleted task forensics:
        - If task XML file is deleted, the registry cache
          may still contain the task GUID and action data
        - $MFT (NTFS master file table) retains metadata
          for deleted task XML files including timestamps
          and original filename
    ossdetect: |-
        Sigma:
        - win_security_scheduled_task_creation.yml
          (parse TaskContent for suspicious actions)
        - proc_creation_win_svchost_susp_child_process.yml
          (svchost Schedule spawning suspicious children)

        Atomic Red Team:
        - T1053.005 (multiple task persistence tests)
        - T1053.005 Test #6 (task pointing to LOLBin)

        Hayabusa:
        - SuspiciousScheduledTaskAction rules
        - TaskSchedulerSvchost child detection

        Velociraptor:
        - Windows.System.ScheduledTasks
          (best single artifact - full task enumeration
          with action path, trigger, and run-as context)
        - Windows.Forensics.Autoruns

        Sysinternals autoruns.exe:
        - Most effective single tool for manual triage
        - Highlights unsigned task actions
        - Shows hidden/broken tasks that PowerShell misses
    tools: |-
        This indicator focuses on HUNTING existing tasks
        rather than catching creation in real-time.

        Useful for:
        - Post-compromise triage (what persists on this host?)
        - Baseline deviation (new tasks since last check)
        - Threat hunting sweeps across fleet

        Tools for task enumeration:
        - Get-ScheduledTask (PowerShell - built-in)
        - schtasks /query /fo LIST /v (schtasks.exe)
        - autoruns.exe -a * -ct (Sysinternals - CSV output)
        - Velociraptor Windows.System.ScheduledTasks
        - Carbon Black / CrowdStrike scheduled task queries

        Red team tools that create tasks:
        SharPersist (/t schtask)
        Cobalt Strike persistence module
        Metasploit persistence/windows/schtasks
        Custom .NET / C++ task COM API callers
          (bypass schtasks.exe entirely)
    actors:
        - APT29
        - APT41
        - Lazarus
    analyst_notes: This indicator is hunting-oriented rather than real-time detection - it's about finding what already exists rather than catching creation live. The two previous T1053.005 indicators cover real-time creation via schtasks.exe EID 1 and Security EID 4698. This one covers the scenario where a task was created by a method that didn't generate obvious telemetry (COM API, fileless stager, pre-Sysmon deployment) and you need to sweep for it. The Get-ScheduledTask PowerShell script and the raw XML scan of C:\Windows\System32\Tasks\ are the most useful live-host forensic tools for this. The svchost.exe spawning suspicious children pattern bridges the gap - it fires when an existing task actually executes, regardless of how it was created, because all modern scheduled tasks run under the Task Scheduler svchost instance. The parent command line containing 'Schedule' distinguishes the Task Scheduler svchost from the dozens of other svchost instances running on a Windows system. Adversaries frequently disguise task names as Windows Update, Windows Defender, Google Update, or similar legitimate-sounding names - always inspect the action path, not the task name, as the primary detection signal.
