# sigma/rules/host/execution/t1053_005_suspicious_task_creation_via_schtasks_exe_supple.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/execution.js
title: 'T1053.005 - Suspicious Task Creation via schtasks.exe (Supplementary: Windows Security Event 4698 (all task creation))'
id: 57d9eba9-facc-52df-9fb9-1c6a27b9fc2d
status: experimental
description: 'schtasks.exe /create with action pointing to suspicious binary path or scripting interpreter Query block: Supplementary: Windows Security Event 4698 (all task creation).'
references:
    - MITRE ATT&CK T1053.005
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.execution
    - attack.t1053.005
    - attack.g0016
    - attack.g0096
    - attack.g0032
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        winlog.event_id: 4698
        winlog.channel: Security
    condition: selection
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/execution.js
    domain: host
    technique: T1053.005
    technique_name: 'Scheduled Task/Job: Scheduled Task'
    os: win
    query_block: 2 of 2
    original_query: |-
        winlog.event_id: 4698
        AND winlog.channel: "Security"
    sysmon: |-
        EventID=1
        Image=*\schtasks.exe
        CommandLine=*/create*
        AND CommandLine matches (any of):
          *\AppData\*
          OR *\Temp\*
          OR *\ProgramData\*
          OR *\Users\Public\*
          OR *powershell*
          OR *cmd.exe*
          OR *wscript*
          OR *cscript*
          OR *mshta*
          OR *rundll32*
          OR *regsvr32*
          OR *certutil*

        // Supplementary - Windows Security Event:
        EventID=4698 (Task Scheduler / Security log)
        // Fires on every task creation regardless of method
        // (schtasks.exe, COM API, remote) - higher coverage
        // than Sysmon EID 1 alone
    powershell: |-
        # Hunt for suspicious schtasks /create executions (Sysmon EID 1)
        Get-WinEvent -FilterHashtable @{
          LogName='Microsoft-Windows-Sysmon/Operational';
          ID=1
        } | Where-Object {
          $_.Properties[4].Value -like '*\schtasks.exe' -and
          $_.Properties[10].Value -match '/create' -and
          $_.Properties[10].Value -match
            '(AppData|Temp|ProgramData|Public|powershell|wscript|cscript|mshta|rundll32|regsvr32|certutil)'
        } | Select TimeCreated,
          @{n='CmdLine';e={$_.Properties[10].Value}},
          @{n='Parent';e={($_.Properties[20].Value -split '\\')[-1]}},
          @{n='User';e={$_.Properties[12].Value}} |
          Sort-Object TimeCreated -Descending

        # Supplementary: Security Event 4698 (task created - all methods)
        Get-WinEvent -FilterHashtable @{
          LogName='Security';
          ID=4698
        } | Select TimeCreated,
          @{n='TaskName';e={
            ([xml]$_.ToXml()).Event.EventData.Data |
            Where-Object { $_.Name -eq 'TaskName' } | Select-Object -Expand '#text'
          }},
          @{n='TaskContent';e={
            ([xml]$_.ToXml()).Event.EventData.Data |
            Where-Object { $_.Name -eq 'TaskContent' } | Select-Object -Expand '#text'
          }},
          @{n='User';e={$_.UserId}} |
          Sort-Object TimeCreated -Descending
    registry: |-
        Scheduled task artifacts on disk:
        - C:\Windows\System32\Tasks\<TaskName>
          (XML file - human readable, contains full task config
          including action, trigger, principal, and run-as user)
        - C:\Windows\SysWOW64\Tasks\<TaskName>
          (32-bit task storage on 64-bit systems)

        Registry task keys (legacy/additional storage):
        - HKLM\SOFTWARE\Microsoft\Windows NT\
          CurrentVersion\Schedule\TaskCache\Tasks\
        - HKLM\SOFTWARE\Microsoft\Windows NT\
          CurrentVersion\Schedule\TaskCache\Tree\
          (tree view of task names by folder path)

        Investigation pivots:
        - Read the XML in C:\Windows\System32\Tasks\
          to see the full task definition:
          Get-Content "C:\Windows\System32\Tasks\<name>"
        - Check SubjectUserName in Security EID 4698 -
          SYSTEM creating tasks is normal;
          a standard user account creating tasks is suspicious
        - Check trigger type: OnLogon and AtStartup triggers
          indicate persistence intent vs one-time execution
    ossdetect: |-
        Sigma:
        - proc_creation_win_schtasks_creation_susp_path.yml
        - proc_creation_win_schtasks_creation_susp_action.yml
        - win_security_scheduled_task_creation.yml (EID 4698)
        - proc_creation_win_schtasks_susp_parent.yml

        Atomic Red Team:
        - T1053.005 Test #1 (schtasks /create local)
        - T1053.005 Test #2 (PowerShell Register-ScheduledTask)
        - Multiple variants covering different triggers/actions

        Hayabusa:
        - ScheduledTaskCreation rules (EID 1 + EID 4698)
        - SuspiciousTaskAction detection category

        Velociraptor:
        - Windows.System.ScheduledTasks
          (enumerates all tasks with full XML content)
        - Windows.EventLogs.Sysmon
        - Windows.Forensics.Autoruns (lists scheduled tasks)

        Sysinternals autoruns.exe:
        - Scheduled Tasks tab shows all registered tasks
        - Highlights unsigned task actions in yellow/red
    tools: |-
        schtasks.exe /create (built-in - most common)
        Task Scheduler MMC (taskschd.msc - GUI, less common in ops)
        PowerShell New-ScheduledTask / Register-ScheduledTask
          (COM-based - does not invoke schtasks.exe,
          only detectable via EID 4698 not EID 1)
        Impacket atexec.py (remote task scheduling)
        CrackMapExec (--exec-method atexec)
        Cobalt Strike (scheduled task persistence module)
        SharPersist (C# persistence tool with schtasks support)

        Common adversary task action patterns:
        - /tr "powershell.exe -ep bypass -w hidden -c <cmd>"
        - /tr "cmd.exe /c <staging command>"
        - /tr "wscript.exe %APPDATA%\payload.vbs"
        - /tr "C:\ProgramData\<random>.exe"
        - /sc ONLOGON /tn "Windows Update" (disguised name)
        - /sc MINUTE /mo 5 (beacon persistence)
    actors:
        - APT29
        - APT41
        - Lazarus
    analyst_notes: 'Scheduled tasks are one of the most commonly abused persistence mechanisms because they are well-understood by IT staff (reducing scrutiny), support a wide range of triggers (logon, startup, time interval, event-based), and can run as SYSTEM with minimal configuration. The schtasks.exe /create command is the most visible path - it generates Sysmon EID 1 with the full command line including the task action, which is often where the suspicious binary path or interpreter is visible. The detection is path-and-action focused: system binaries in standard paths (C:\Windows\System32) creating tasks that run system management tools are normal; tasks running interpreters (powershell, wscript, mshta) or binaries from user-writable paths (AppData, Temp, ProgramData) are suspicious. Security Event 4698 is the more comprehensive signal because it fires regardless of how the task was created - schtasks.exe, PowerShell COM API, or remote scheduling - and the event content includes the full task XML. The task XML file in C:\Windows\System32\Tasks\ is a valuable forensic artifact: it contains the full action, trigger, principal, and run-as context and persists even after the process that created it has exited. Adversaries frequently disguise task names as legitimate Windows tasks (''Windows Update'', ''GoogleUpdateTask'', ''AdobeFlashUpdate'') - the action path is more reliable than the task name for detection.'
