# sigma/rules/host/execution/t1569_002_service_binary_in_non_standard_path_real_time_im.yml
# Generated by TONK sigma-export.py - do not edit by hand.
# Source of truth: host/js/data/execution.js
title: 'T1569.002 - Service Binary in Non-Standard Path (Real-time: ImagePath written outside standard paths)'
id: 27863851-4b97-511e-88e2-2a028d07dd99
status: experimental
description: 'Registered service with ImagePath outside standard system directories - hunting for malicious persistence via service registry Query block: Real-time: ImagePath written outside standard paths.'
references:
    - MITRE ATT&CK T1569.002
    - T1543.003
author: TONK Hunt Reference (6b74.dev)
date: 2026/07/21
tags:
    - attack.execution
    - attack.t1569.002
    - attack.g0016
    - attack.g0096
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        winlog.event_id: 13
        registry.path: '*\Services\*\ImagePath'
    filter_0:
        registry.data.strings:
            - C:\\Windows\\System32\\*
            - C:\\Windows\\SysWOW64\\*
            - C:\\Program Files\\*
            - C:\\Program Files (x86)\\*
    condition: selection and not filter_0
falsepositives:
    - Unknown - baseline before enabling
level: medium
tonk:
    fidelity: full
    source_file: host/js/data/execution.js
    domain: host
    technique: T1569.002
    technique_name: 'System Services: Service Execution'
    os: win
    query_block: 1 of 2
    original_query: |-
        winlog.event_id: 13
        AND registry.path: *\Services\*\ImagePath
        AND NOT registry.data.strings: ("C:\\Windows\\System32\\*" OR "C:\\Windows\\SysWOW64\\*" OR "C:\\Program Files\\*" OR "C:\\Program Files (x86)\\*")
    sysmon: |-
        // Real-time: Sysmon EID 13 (registry value set)
        // catches service ImagePath being written:
        EventID=13
        TargetObject=*\Services\*\ImagePath
        Details NOT matching:
          C:\Windows\System32\*
          C:\Windows\SysWOW64\*
          C:\Program Files\*
          C:\Program Files (x86)\*

        // Supplementary - service execution artifact:
        EventID=1
        ParentImage=*\services.exe
        Image NOT matching known-good service paths
        // services.exe is the SCM host - it spawns
        // services directly, making it a useful parent filter
    powershell: |-
        # Hunt for services with ImagePath outside standard directories
        $standardPaths = @(
          'C:\Windows\',
          'C:\Program Files\',
          'C:\Program Files (x86)\'
        )

        Get-ItemProperty "HKLM:\SYSTEM\CurrentControlSet\Services\*" |
          Where-Object {
            $ip = $_.ImagePath
            if (-not $ip) { return $false }
            $ip = $ip.TrimStart('"') -replace '^([^"]+).*','$1'
            $ip = [System.Environment]::ExpandEnvironmentVariables($ip)
            -not ($standardPaths | Where-Object { $ip -like "$_*" })
          } | Select-Object `
            @{n='ServiceName';e={$_.PSChildName}},
            @{n='ImagePath';e={$_.ImagePath}},
            @{n='Start';e={
              switch ($_.Start) {
                0 {'Boot'} 1 {'System'} 2 {'Auto'}
                3 {'Demand'} 4 {'Disabled'} default {$_}
              }
            }},
            @{n='ObjectName';e={$_.ObjectName}},
            @{n='DisplayName';e={$_.DisplayName}} |
          Sort-Object ServiceName
    registry: |-
        Primary artifact - service ImagePath:
        HKLM\SYSTEM\CurrentControlSet\Services\<Name>\
          ImagePath - full path to service binary
          (may include arguments after the binary path)
          (may be quoted or unquoted)
          (may use environment variables: %SystemRoot%\...)

        Suspicious ImagePath patterns to hunt:
        - C:\Users\<user>\AppData\*
        - C:\Windows\Temp\*
        - C:\ProgramData\*
        - C:\Users\Public\*
        - Relative paths (no drive letter)
        - UNC paths (\\server\share\binary.exe)
        - Paths containing scripting interpreters as the binary

        Baseline approach:
        - Export all service ImagePaths on a known-good system
        - Compare against target system
        - New entries with non-standard paths = investigate

        Deleted service forensics:
        - System Event 7036 (service stopped) and
          absence of 7045 for a service that previously
          existed suggests manual deletion
        - MFT and USN Journal retain file metadata
          for the service binary even after deletion
    ossdetect: |-
        Sigma:
        - registry_event_win_service_registry_susp_path.yml
        - proc_creation_win_services_susp_child_process.yml

        Atomic Red Team:
        - T1569.002 (service creation variants)
        - T1543.003 (Create or Modify System Process: Windows Service)

        Hayabusa:
        - ServiceBinaryNonStandardPath rules (EID 13)
        - services.exe child process anomaly detection

        Velociraptor:
        - Windows.System.Services
          (best single artifact for fleet-wide service hunting)
        - Windows.Forensics.Autoruns

        Sysinternals autoruns.exe:
        - Services tab with VirusTotal integration
        - Highlights unsigned service binaries in red/yellow
        - Most practical tool for manual IR triage
    tools: |-
        This indicator is hunting-oriented - sweep for
        existing malicious services rather than real-time
        creation detection (covered in indicator 1).

        Most useful in:
        - IR triage: what services exist that shouldn't?
        - Baseline deviation: new services since last audit
        - Post-persistence-establishment hunting
        - Fleet-wide sweep via Velociraptor or EDR

        Service creation tools (all leave ImagePath artifact):
        sc.exe create
        PowerShell New-Service / Register-ServiceJob
        Compiled code using CreateService() Win32 API
        Metasploit persistence modules
        SharPersist
        Cobalt Strike persistence
        Any tool writing directly to the Services registry key

        Known-good exceptions to tune out:
        - Security software (AV, EDR agents)
        - Monitoring agents (Elastic, Splunk UF, etc.)
        - Vendor software with non-standard install paths
        - Build allowlist from asset management / SCCM data
    actors:
        - APT29
        - APT41
        - FIN6
    analyst_notes: 'This is the hunt-oriented companion to the real-time sc.exe detection - same concept as the scheduled task action hunting indicator in T1053.005. The PowerShell registry sweep is the core tool: it enumerates every registered service''s ImagePath and flags anything outside standard system and program directories. The key challenge is normalization: ImagePath values can be quoted, unquoted, contain environment variables (%SystemRoot%), or embed arguments after the binary path - the hunt script handles this. Environment variable expansion is important because %SystemRoot%\system32\svchost.exe is legitimate but looks non-standard if you string-match on C:\Windows without expanding the variable first. Services running from C:\ProgramData\ deserve particular scrutiny - it''s a world-writable directory that doesn''t appear in standard path allowlists, making it a common adversary staging location. The services.exe parent filter in the Sysmon section is a useful real-time complement: services.exe is the SCM binary that directly spawns services, so filtering Sysmon EID 1 on ParentImage=services.exe and checking the child''s path against standard directories gives live execution coverage for this indicator.'
